# cvs/src/etc/rc.d/iked,v
head	1.3;
access;
symbols
	OPENBSD_6_1:1.3.0.8
	OPENBSD_6_1_BASE:1.3
	OPENBSD_6_0:1.3.0.6
	OPENBSD_6_0_BASE:1.3
	OPENBSD_5_9:1.3.0.2
	OPENBSD_5_9_BASE:1.3
	OPENBSD_5_8:1.1.0.20
	OPENBSD_5_8_BASE:1.1
	OPENBSD_5_7:1.1.0.12
	OPENBSD_5_7_BASE:1.1
	OPENBSD_5_6:1.1.0.16
	OPENBSD_5_6_BASE:1.1
	OPENBSD_5_5:1.1.0.14
	OPENBSD_5_5_BASE:1.1
	OPENBSD_5_4:1.1.0.10
	OPENBSD_5_4_BASE:1.1
	OPENBSD_5_3:1.1.0.8
	OPENBSD_5_3_BASE:1.1
	OPENBSD_5_2:1.1.0.6
	OPENBSD_5_2_BASE:1.1
	OPENBSD_5_1_BASE:1.1
	OPENBSD_5_1:1.1.0.4
	OPENBSD_5_0:1.1.0.2
	OPENBSD_5_0_BASE:1.1;
locks; strict;
comment	@# @;


1.3
date	2015.12.19.13.45.12;	author ajacoutot;	state Exp;
branches;
next	1.2;
commitid	Fo7f7cSN2lRe42tW;

1.2
date	2015.12.19.10.09.04;	author ajacoutot;	state Exp;
branches;
next	1.1;
commitid	gZyqXQTxi4wOLgzw;

1.1
date	2011.07.06.18.55.36;	author robert;	state Exp;
branches;
next	;


desc
@@


1.3
log
@Child will not return a config parsing error to the parent so rc.d will lie and
tell you that the daemon succesfully started; prevent this by running a config
check first.
Note that we only do this for those particular daemons, not all that support a
config test mode.

ok benno@@ phessler@@
@
text
@#!/bin/sh
#
# $OpenBSD: iked,v 1.2 2015/12/19 10:09:04 ajacoutot Exp $

daemon="/sbin/iked"

. /etc/rc.d/rc.subr

pexp="iked: parent.*"

rc_pre() {
	[ X"${sasyncd_flags}" != X"NO" ] && \
		daemon_flags="-S ${daemon_flags}"
	#return 0
	# child will not return a config parsing error to the parent
	${daemon} -n ${daemon_flags}
}

rc_cmd $1
@


1.2
log
@Tweak pexp.
@
text
@d3 1
a3 1
# $OpenBSD: iked,v 1.1 2011/07/06 18:55:36 robert Exp $
d14 3
a16 1
	return 0
@


1.1
log
@Add rc.d(8) script for the system daemons that are restartable.
From now on rc(8) is going to call these scripts to start them up on boot
in the same order than before.
In addition the inetd and rwhod variables in rc.conf are deprecated so that
inetd_flags and rwhod_flags should be used. The old flags are still going
to be used for some time to allow users to switch.
There are more rc modifications to come later so let's put this in so
we can base more work on this.
It is important to mention that you can still keep using rc.local just
like the way you did before, and we have no intention to remove that either.

I'd also like to thank ajacoutot@@, halex@@, sthen@@ and schwarze@@ for working
on this with me.
@
text
@d3 1
a3 1
# $OpenBSD$
d9 1
a9 1
pexp="iked:.*parent.*"
@

