# cvs/src/etc/rc.d/ntpd,v
head	1.3;
access;
symbols
	OPENBSD_6_2:1.3.0.10
	OPENBSD_6_2_BASE:1.3
	OPENBSD_6_1:1.3.0.8
	OPENBSD_6_1_BASE:1.3
	OPENBSD_6_0:1.3.0.6
	OPENBSD_6_0_BASE:1.3
	OPENBSD_5_9:1.3.0.2
	OPENBSD_5_9_BASE:1.3
	OPENBSD_5_8:1.2.0.18
	OPENBSD_5_8_BASE:1.2
	OPENBSD_5_7:1.2.0.10
	OPENBSD_5_7_BASE:1.2
	OPENBSD_5_6:1.2.0.14
	OPENBSD_5_6_BASE:1.2
	OPENBSD_5_5:1.2.0.12
	OPENBSD_5_5_BASE:1.2
	OPENBSD_5_4:1.2.0.8
	OPENBSD_5_4_BASE:1.2
	OPENBSD_5_3:1.2.0.6
	OPENBSD_5_3_BASE:1.2
	OPENBSD_5_2:1.2.0.4
	OPENBSD_5_2_BASE:1.2
	OPENBSD_5_1_BASE:1.2
	OPENBSD_5_1:1.2.0.2
	OPENBSD_5_0:1.1.0.2
	OPENBSD_5_0_BASE:1.1;
locks; strict;
comment	@# @;


1.3
date	2016.02.02.17.51.11;	author sthen;	state Exp;
branches;
next	1.2;
commitid	LQHdRRMF0nq0fugE;

1.2
date	2011.09.14.02.36.09;	author william;	state Exp;
branches;
next	1.1;

1.1
date	2011.07.06.18.55.36;	author robert;	state Exp;
branches;
next	;


desc
@@


1.3
log
@Remove setproctitle() for the parent process. Because rc.d(8) uses process
titles (including flags) to distinguish between daemons, this makes it
possible to manage multiple copies of a daemon using the normal infrastructure
by symlinking rc.d scripts to a new name. ok jung@@ ajacoutot@@, smtpd ok gilles@@
@
text
@#!/bin/sh
#
# $OpenBSD: ntpd,v 1.2 2011/09/14 02:36:09 william Exp $

daemon="/usr/sbin/ntpd"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
@


1.2
log
@ntpd doesn't really support reload, so don't allow it in the rc script.

ok ajacoutot robert
@
text
@d3 1
a3 1
# $OpenBSD: ntpd,v 1.1 2011/07/06 18:55:36 robert Exp $
a8 1
pexp="ntpd: \[priv\]"
@


1.1
log
@Add rc.d(8) script for the system daemons that are restartable.
From now on rc(8) is going to call these scripts to start them up on boot
in the same order than before.
In addition the inetd and rwhod variables in rc.conf are deprecated so that
inetd_flags and rwhod_flags should be used. The old flags are still going
to be used for some time to allow users to switch.
There are more rc modifications to come later so let's put this in so
we can base more work on this.
It is important to mention that you can still keep using rc.local just
like the way you did before, and we have no intention to remove that either.

I'd also like to thank ajacoutot@@, halex@@, sthen@@ and schwarze@@ for working
on this with me.
@
text
@d3 1
a3 1
# $OpenBSD$
d10 1
@

