# cvs/src/etc/rc.d/slaacd,v
head	1.2;
access;
symbols;
locks; strict;
comment	@# @;


1.2
date	2017.06.03.10.00.29;	author florian;	state Exp;
branches;
next	1.1;
commitid	LiVl5Mw0RdD7hl8R;

1.1
date	2017.05.29.10.24.06;	author florian;	state Exp;
branches;
next	;
commitid	GMPrMXlLo65fu9RV;


desc
@@


1.2
log
@Move slaacd to /sbin
jca points out that all the other interface configuration tools live
there (like ifconfig or dhclient). Furthermore it starts so early in
the boot process that /usr might not be mounted yet if it's a nfs
filesystem.
sthen and deraadt agree
@
text
@#!/bin/sh
#
# $OpenBSD: slaacd,v 1.1 2017/05/29 10:24:06 florian Exp $

daemon="/sbin/slaacd"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
@


1.1
log
@rc.d(8) for slaacd
OK phessler, deraadt
@
text
@d3 1
a3 1
# $OpenBSD$
d5 1
a5 1
daemon="/usr/sbin/slaacd"
@

