# cvs/src/etc/rc.d/slaacd,v
head	1.1;
access;
symbols;
locks; strict;
comment	@# @;


1.1
date	2017.05.29.10.24.06;	author florian;	state Exp;
branches;
next	;
commitid	GMPrMXlLo65fu9RV;


desc
@@


1.1
log
@rc.d(8) for slaacd
OK phessler, deraadt
@
text
@#!/bin/sh
#
# $OpenBSD$

daemon="/usr/sbin/slaacd"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
@
