# cvs/src/etc/rc.d/tftpproxy,v
head	1.2;
access;
symbols
	OPENBSD_6_0:1.2.0.8
	OPENBSD_6_0_BASE:1.2
	OPENBSD_5_9:1.2.0.2
	OPENBSD_5_9_BASE:1.2
	OPENBSD_5_8:1.2.0.4
	OPENBSD_5_8_BASE:1.2
	OPENBSD_5_7:1.1.0.6
	OPENBSD_5_7_BASE:1.1
	OPENBSD_5_6:1.1.0.10
	OPENBSD_5_6_BASE:1.1
	OPENBSD_5_5:1.1.0.8
	OPENBSD_5_5_BASE:1.1
	OPENBSD_5_4:1.1.0.4
	OPENBSD_5_4_BASE:1.1
	OPENBSD_5_3:1.1.0.2
	OPENBSD_5_3_BASE:1.1;
locks; strict;
comment	@# @;


1.2
date	2015.07.15.13.50.54;	author ajacoutot;	state Exp;
branches;
next	1.1;
commitid	oeNl5gMltmsYwI9M;

1.1
date	2012.08.20.12.37.53;	author dlg;	state Exp;
branches;
next	;


desc
@@


1.2
log
@Fix base rc.d scripts after the recent rc.subr change.

ok halex@@
@
text
@#!/bin/sh
#
# $OpenBSD: tftpproxy,v 1.1 2012/08/20 12:37:53 dlg Exp $

daemon="/usr/sbin/tftp-proxy"

. /etc/rc.d/rc.subr

pexp="tftp-proxy: privproc"
rc_reload=NO

rc_cmd $1
@


1.1
log
@wire up the bits for tftp-proxy, which is called tftpproxy from rc.d's
point of view.

mostly ok sthen@@ ajacoutot@@, who were discussing the feng shui of the
start_daemons chunk
@
text
@d3 1
a3 1
# $OpenBSD$
d9 1
@

