# cvs/src/etc/rc.d/vmd,v
head	1.4;
access;
symbols
	OPENBSD_6_2:1.4.0.10
	OPENBSD_6_2_BASE:1.4
	OPENBSD_6_1:1.4.0.8
	OPENBSD_6_1_BASE:1.4
	OPENBSD_6_0:1.4.0.6
	OPENBSD_6_0_BASE:1.4
	OPENBSD_5_9:1.4.0.2
	OPENBSD_5_9_BASE:1.4;
locks; strict;
comment	@# @;


1.4
date	2016.02.02.17.51.11;	author sthen;	state Exp;
branches;
next	1.3;
commitid	LQHdRRMF0nq0fugE;

1.3
date	2015.12.19.13.45.12;	author ajacoutot;	state Exp;
branches;
next	1.2;
commitid	Fo7f7cSN2lRe42tW;

1.2
date	2015.12.16.14.21.28;	author jturner;	state Exp;
branches;
next	1.1;
commitid	3n4ZaaXNXoODlJw6;

1.1
date	2015.12.05.09.55.18;	author jasper;	state Exp;
branches;
next	;
commitid	BDIjpYbNGByke3hj;


desc
@@


1.4
log
@Remove setproctitle() for the parent process. Because rc.d(8) uses process
titles (including flags) to distinguish between daemons, this makes it
possible to manage multiple copies of a daemon using the normal infrastructure
by symlinking rc.d scripts to a new name. ok jung@@ ajacoutot@@, smtpd ok gilles@@
@
text
@#!/bin/sh
#
# $OpenBSD: vmd,v 1.3 2015/12/19 13:45:12 ajacoutot Exp $

daemon="/usr/sbin/vmd"

. /etc/rc.d/rc.subr

# child will not return a config parsing error to the parent
rc_pre() {
	${daemon} -n ${daemon_flags}
}

rc_cmd $1
@


1.3
log
@Child will not return a config parsing error to the parent so rc.d will lie and
tell you that the daemon succesfully started; prevent this by running a config
check first.
Note that we only do this for those particular daemons, not all that support a
config test mode.

ok benno@@ phessler@@
@
text
@d3 1
a3 1
# $OpenBSD: vmd,v 1.2 2015/12/16 14:21:28 jturner Exp $
a7 2

pexp="vmd: parent.*"
@


1.2
log
@Add pexp match otherwise rcctl stop vmd doesn't work. OK deraadt@@
@
text
@d3 1
a3 1
# $OpenBSD: vmd,v 1.1 2015/12/05 09:55:18 jasper Exp $
d10 5
@


1.1
log
@add rcscript for vmd

ok mlarkin@@
prompted by deraadt@@
@
text
@d3 1
a3 1
# $OpenBSD$
d8 2
@

