# cvs/ports/lang/php/5.4/pkg/Attic/php54_fpm.rc,v
head	1.2;
access;
symbols
	OPENBSD_5_9:1.1.0.2
	OPENBSD_5_9_BASE:1.1;
locks; strict;
comment	@ * @;


1.2
date	2016.04.28.18.06.47;	author sthen;	state dead;
branches;
next	1.1;
commitid	i47ZnKZZ6MYLicYj;

1.1
date	2015.09.18.09.18.22;	author sthen;	state Exp;
branches;
next	;
commitid	x1OqzrB6vMfHx6Ez;


desc
@@


1.2
log
@5.4
@
text
@#!/bin/sh
#
# $OpenBSD: php54_fpm.rc,v 1.1 2015/09/18 09:18:22 sthen Exp $

daemon="${TRUEPREFIX}/sbin/php-fpm-${PV}"

. /etc/rc.d/rc.subr

pexp="php-fpm-${PV}: master process .*"
rc_reload=NO

rc_cmd $1
@


1.1
log
@de-conflict php-fpm packages, ok robert@@

Heads up! The rc script has changed name:

  sed -i s/php_fpm/php56_fpm/g /etc/rc.conf.local
@
text
@d3 1
a3 1
# $OpenBSD: php_fpm.rc,v 1.2 2015/07/15 16:03:07 ajacoutot Exp $
@

