# cvs/ports/sysutils/anacron/pkg/Attic/anacron.rc,v
head	1.2;
access;
symbols
	OPENBSD_5_8:1.1.0.14
	OPENBSD_5_8_BASE:1.1
	OPENBSD_5_7:1.1.0.10
	OPENBSD_5_7_BASE:1.1
	OPENBSD_5_6:1.1.0.8
	OPENBSD_5_6_BASE:1.1
	OPENBSD_5_5:1.1.0.6
	OPENBSD_5_5_BASE:1.1
	OPENBSD_5_4:1.1.0.4
	OPENBSD_5_4_BASE:1.1
	OPENBSD_5_3:1.1.0.2
	OPENBSD_5_3_BASE:1.1;
locks; strict;
comment	@ * @;


1.2
date	2016.01.11.06.59.48;	author ajacoutot;	state dead;
branches;
next	1.1;
commitid	um0J6dSULjpjAoqE;

1.1
date	2012.11.20.15.14.41;	author giovanni;	state Exp;
branches;
next	;


desc
@@


1.2
log
@anacron(8) is not a daemon, so there's no reason to have an rc.d script.
Instead use @@reboot in root's crontab.

ok giovanni@@ (maintainer)
@
text
@#!/bin/sh
#
# $OpenBSD: anacron.rc,v 1.1 2012/11/20 15:14:41 giovanni Exp $

daemon="${TRUEPREFIX}/sbin/anacron"
daemon_flags="-s"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
@


1.1
log
@
Add an rcscript and move MESSAGE to README
take maintainership
ok Andreas Kahari (old maintainer), gonzalo@@
@
text
@d3 1
a3 1
# $OpenBSD$
@

