# cvs/ports/devel/go-tools/pkg/godoc.rc,v
head	1.1;
branch	1.1.1;
access;
symbols
	OPENBSD_6_1:1.1.1.1.0.8
	OPENBSD_6_1_BASE:1.1.1.1
	OPENBSD_6_0:1.1.1.1.0.6
	OPENBSD_6_0_BASE:1.1.1.1
	OPENBSD_5_9:1.1.1.1.0.2
	OPENBSD_5_9_BASE:1.1.1.1
	OPENBSD_5_8:1.1.1.1.0.4
	OPENBSD_5_8_BASE:1.1.1.1
	czarkoff_20150608:1.1.1.1
	czarkoff:1.1.1;
locks; strict;
comment	@ * @;


1.1
date	2015.06.08.09.41.03;	author czarkoff;	state Exp;
branches
	1.1.1.1;
next	;
commitid	pWx0CLB6XHkGDmny;

1.1.1.1
date	2015.06.08.09.41.03;	author czarkoff;	state Exp;
branches;
next	;
commitid	pWx0CLB6XHkGDmny;


desc
@@



1.1
log
@Initial revision
@
text
@#!/bin/sh
#
# $OpenBSD: godoc.rc,v 1.2 2015/03/21 19:44:36 ajacoutot Exp $

daemon="${TRUEPREFIX}/bin/godoc"
daemon_flags="-http=localhost:6060 -v=true"
daemon_user=www

. /etc/rc.d/rc.subr

rc_bg=YES
rc_reload=NO

rc_cmd $1
@


1.1.1.1
log
@import go-tools

This package holds the source for various packages and tools that support the Go
programming language:

 * benchcmp	displays performance changes between benchmarks
 * callgraph	reports the call graph of a Go program
 * digraph	performs queries over unlabelled directed graphs represented in
		text form
 * eg		performs example-based refactoring
 * fiximports	fixes import declarations to use the canonical import path for
 		packages that have an "import comment"
 * godex	dumps exported information of packages or selected objects
 * godoc	extracts and generates documentation for Go programs
 * goimports	updates Go import lines
 * gomvpkg	moves go packages, updating import declarations
 * gorename	performs precise type-safe renaming of identifiers in Go code
 * gotype	does syntactic and semantic analysis of Go files and packages
 * oracle	answers questions about Go source code
 * ssadump	displays and interprets the SSA form of Go programs
 * stress	runs a given process in parallel in a loop and collects failures
 * stringer	creats of methods that satisfy the fmt.Stringer interface
 * cover	analyzes the coverage profiles
 * vet		examines Go source code and reports suspicious constructs

OK sthen@@
@
text
@@
