# cvs/ports/net/ices2/pkg/ices2.rc,v
head	1.2;
access;
symbols
	OPENBSD_6_1:1.1.0.8
	OPENBSD_6_1_BASE:1.1
	OPENBSD_6_0:1.1.0.6
	OPENBSD_6_0_BASE:1.1
	OPENBSD_5_9:1.1.0.2
	OPENBSD_5_9_BASE:1.1
	OPENBSD_5_8:1.1.0.4
	OPENBSD_5_8_BASE:1.1;
locks; strict;
comment	@ * @;


1.2
date	2017.05.27.13.53.07;	author ajacoutot;	state Exp;
branches;
next	1.1;
commitid	SDpZdO8hfYIlMQ9j;

1.1
date	2015.04.15.10.13.21;	author jasper;	state Exp;
branches;
next	;
commitid	Srqy1efJNcxit1uR;


desc
@@


1.2
log
@Cope with recent rc.subr(8) change (removal of rc_bg etc.)
@
text
@#!/bin/sh
#
# $OpenBSD: ices2.rc,v 1.1 2015/04/15 10:13:21 jasper Exp $

daemon="${TRUEPREFIX}/bin/ices2"
daemon_flags="${SYSCONFDIR}/ices2.xml"

. /etc/rc.d/rc.subr

rc_cmd $1
@


1.1
log
@- update to ices-2.0.2
- reformat Makefile to make it readable again
- use correct AUTOCONF_VERSION
- add rc script

ok aja@@
@
text
@d3 1
a3 1
# $OpenBSD$
a8 2

rc_bg=YES
@

