# cvs/ports/comms/lcdproc/pkg/Attic/LCDd.rc,v
head	1.2;
access;
symbols
	edd_20111201:1.1.1.1 edd:1.1.1;
locks; strict;
comment	@ * @;


1.2
date	2011.12.02.12.39.23;	author edd;	state dead;
branches;
next	1.1;

1.1
date	2011.12.01.19.05.11;	author edd;	state Exp;
branches
	1.1.1.1;
next	;

1.1.1.1
date	2011.12.01.19.05.11;	author edd;	state Exp;
branches;
next	;


desc
@@


1.2
log
@lowercase LCDd rc script filename at request of aja@@
@
text
@#!/bin/sh
#
# $OpenBSD: LCDd.rc,v 1.1 2011/12/01 19:05:11 edd Exp $

daemon="${TRUEPREFIX}/sbin/LCDd"

. /etc/rc.d/rc.subr

rc_cmd $1
@


1.1
log
@Initial revision
@
text
@d3 1
a3 1
# $OpenBSD$
@


1.1.1.1
log
@import lcdproc - a tool for interfacing LCD modules

OK sebastia@@

@
text
@@
