# cvs/ports/www/apache-httpd/pkg/Attic/httpd2.rc,v
head	1.2;
access;
symbols
	OPENBSD_5_5:1.1.0.8
	OPENBSD_5_5_BASE:1.1
	OPENBSD_5_4:1.1.0.6
	OPENBSD_5_4_BASE:1.1
	OPENBSD_5_3:1.1.0.4
	OPENBSD_5_3_BASE:1.1
	OPENBSD_5_2:1.1.0.2
	OPENBSD_5_2_BASE:1.1;
locks; strict;
comment	@ * @;


1.2
date	2014.07.27.10.51.37;	author kili;	state dead;
branches;
next	1.1;
commitid	gYt5lPsCZtEp8W4S;

1.1
date	2012.07.18.08.45.42;	author sthen;	state Exp;
branches;
next	;


desc
@@


1.2
log
@Rename rc scripts of apache-httpd to apache resp. apache2 to get rid
of a conflict with base httpd.

ok sthen@@
@
text
@#!/bin/sh
#
# $OpenBSD: httpd2.rc,v 1.1 2012/07/18 08:45:42 sthen Exp $

daemon="${TRUEPREFIX}/sbin/httpd2"

. /etc/rc.d/rc.subr

rc_pre() {
	rm -f /var/apache2/logs/{ssl_mutex,httpd.lock,accept.lock}.*
}

rc_cmd $1
@


1.1
log
@Update to the Apache HTTP server 2.2.22, the current release on the 2.2 branch.
Merged from diffs from myself, Rodolfo Gouveia and with an rc script from
giovanni@@.

This is a security update, fixing CVE-2011-3348, CVE-2011-3368, CVE-2011-3607,
CVE-2011-4317, CVE-2012-0021, CVE-2012-0031 and CVE-2012-0053.
@
text
@d3 1
a3 1
# $OpenBSD: httpd,v 1.1 2011/07/06 18:55:36 robert Exp $
@

