# cvs/ports/geo/mapserver/pkg/mapserv.rc,v
head	1.2;
access;
symbols
	OPENBSD_6_0:1.2.0.6
	OPENBSD_6_0_BASE:1.2
	OPENBSD_5_9:1.2.0.2
	OPENBSD_5_9_BASE:1.2
	OPENBSD_5_8:1.2.0.4
	OPENBSD_5_8_BASE:1.2
	OPENBSD_5_7:1.1.0.4
	OPENBSD_5_7_BASE:1.1
	OPENBSD_5_6:1.1.0.2
	OPENBSD_5_6_BASE:1.1;
locks; strict;
comment	@ * @;


1.2
date	2015.07.15.23.56.36;	author ajacoutot;	state Exp;
branches;
next	1.1;
commitid	7OcA7Rb7FQtD8Tzh;

1.1
date	2014.05.18.09.29.50;	author landry;	state Exp;
branches;
next	;


desc
@@


1.2
log
@Fix pexp.
@
text
@#!/bin/sh
#
# $OpenBSD: mapserv.rc,v 1.1 2014/05/18 09:29:50 landry Exp $

daemon="${LOCALBASE}/bin/spawn-fcgi"
daemon_flags="-c ${VARBASE}/www -s run/mapserv.sock -u www -- /cgi-bin/mapserv"

. /etc/rc.d/rc.subr

pexp=".*/cgi-bin/mapserv"

rc_cmd $1
@


1.1
log
@Bugfix update to mapserver 6.4.1.

Adapt README's for nginx, and provide an rc script using spawn-fcgi.
@
text
@d3 1
a3 1
# $OpenBSD$
d10 1
a10 1
pexp="/cgi-bin/mapserv"
@

