# cvs/ports/x11/dbus/pkg/Attic/dbus-daemon.rc,v
head	1.2;
access;
symbols;
locks; strict;
comment	@ * @;


1.2
date	2010.10.28.21.10.07;	author robert;	state dead;
branches;
next	1.1;

1.1
date	2010.10.28.10.59.33;	author robert;	state Exp;
branches;
next	;


desc
@@


1.2
log
@After discussion with $french and $german, we are renaming the rc scripts from
dash to underscore so that we are consitent with the foo_bar_flags
variable.
@
text
@#!/bin/sh

. /etc/rc.d/rc.subr

daemon="${TRUEPREFIX}/bin/dbus-daemon"
daemon_flags="--system"

rc_pre() {
	install -d -o _dbus -g _dbus /var/run/dbus
}

rc_post() {
	rm /var/run/dbus/pid
}

rc_cmd $1
@


1.1
log
@add an rc script for dbus-daemon
@
text
@@

