# cvs/ports/security/heimdal/pkg/kpasswdd.rc,v
head	1.3;
access;
symbols
	OPENBSD_6_2:1.3.0.2
	OPENBSD_6_2_BASE:1.3
	OPENBSD_6_1:1.1.0.12
	OPENBSD_6_1_BASE:1.1
	OPENBSD_6_0:1.1.0.10
	OPENBSD_6_0_BASE:1.1
	OPENBSD_5_9:1.1.0.6
	OPENBSD_5_9_BASE:1.1
	OPENBSD_5_8:1.1.0.8
	OPENBSD_5_8_BASE:1.1
	OPENBSD_5_7:1.1.0.4
	OPENBSD_5_7_BASE:1.1
	OPENBSD_5_6:1.1.0.2
	OPENBSD_5_6_BASE:1.1;
locks; strict;
comment	@ * @;


1.3
date	2017.05.27.19.57.06;	author ajacoutot;	state Exp;
branches;
next	1.2;
commitid	H08QxD8pqoInz6Rz;

1.2
date	2017.05.27.13.53.07;	author ajacoutot;	state Exp;
branches;
next	1.1;
commitid	SDpZdO8hfYIlMQ9j;

1.1
date	2014.07.13.14.10.14;	author ajacoutot;	state Exp;
branches;
next	;
commitid	R7CQhUQNdJfZrGZj;


desc
@@


1.3
log
@Bring rc_bg back.
@
text
@#!/bin/sh
#
# $OpenBSD: kpasswdd.rc,v 1.2 2017/05/27 13:53:07 ajacoutot Exp $

daemon="${TRUEPREFIX}/heimdal/libexec/kpasswdd"

. /etc/rc.d/rc.subr

rc_bg=YES
rc_reload=NO

rc_cmd $1
@


1.2
log
@Cope with recent rc.subr(8) change (removal of rc_bg etc.)
@
text
@d3 1
a3 1
# $OpenBSD: kpasswdd.rc,v 1.1 2014/07/13 14:10:14 ajacoutot Exp $
d9 1
@


1.1
log
@Import heimdal-1.5.3; aka resurrect security/heimdal.

Heimdal is an implementation of Kerberos 5 (and some more stuff) largely
written in Sweden.

ok sthen@@

Note that it is not enabled yet because none of the requesters for this
port stood up to give me any feedback nor OK. Also this is still a WIP, it
may not even build with all the recent changes in base.
@
text
@d3 1
a3 1
# $OpenBSD: kpasswdd,v 1.1 2011/07/07 18:17:50 ajacoutot Exp $
a8 1
rc_bg=YES
@

