# cvs/ports/www/cherokee/pkg/cherokee.rc,v
head	1.6;
access;
symbols
	OPENBSD_6_0:1.6.0.6
	OPENBSD_6_0_BASE:1.6
	OPENBSD_5_9:1.6.0.2
	OPENBSD_5_9_BASE:1.6
	OPENBSD_5_8:1.6.0.4
	OPENBSD_5_8_BASE:1.6
	OPENBSD_5_7:1.5.0.10
	OPENBSD_5_7_BASE:1.5
	OPENBSD_5_6:1.5.0.8
	OPENBSD_5_6_BASE:1.5
	OPENBSD_5_5:1.5.0.6
	OPENBSD_5_5_BASE:1.5
	OPENBSD_5_4:1.5.0.4
	OPENBSD_5_4_BASE:1.5
	OPENBSD_5_3:1.5.0.2
	OPENBSD_5_3_BASE:1.5
	OPENBSD_5_2:1.4.0.4
	OPENBSD_5_2_BASE:1.4
	OPENBSD_5_1_BASE:1.4
	OPENBSD_5_1:1.4.0.2
	OPENBSD_5_0:1.3.0.4
	OPENBSD_5_0_BASE:1.3
	OPENBSD_4_9:1.3.0.2
	OPENBSD_4_9_BASE:1.3;
locks; strict;
comment	@ * @;


1.6
date	2015.07.15.23.29.57;	author ajacoutot;	state Exp;
branches;
next	1.5;
commitid	vq4hx7arUqZcQCc0;

1.5
date	2012.12.12.02.17.40;	author brad;	state Exp;
branches;
next	1.4;

1.4
date	2011.12.17.22.01.13;	author sthen;	state Exp;
branches;
next	1.3;

1.3
date	2010.12.27.14.50.24;	author ajacoutot;	state Exp;
branches;
next	1.2;

1.2
date	2010.12.24.10.40.06;	author ajacoutot;	state Exp;
branches;
next	1.1;

1.1
date	2010.12.20.07.51.58;	author ajacoutot;	state Exp;
branches;
next	;


desc
@@


1.6
log
@Fix daemon and pexp.
@
text
@#!/bin/sh
#
# $OpenBSD: cherokee.rc,v 1.5 2012/12/12 02:17:40 brad Exp $

daemon="${TRUEPREFIX}/sbin/cherokee -d"

. /etc/rc.d/rc.subr

# extra space injected in the process list
pexp="${TRUEPREFIX}/sbin/cherokee ${daemon_flags:+ ${daemon_flags}}"

rc_cmd $1
@


1.5
log
@Upgrade to cherokee 1.2.101.

ok sthen@@
@
text
@d3 1
a3 1
# $OpenBSD: cherokee.rc,v 1.4 2011/12/17 22:01:13 sthen Exp $
d5 1
a5 2
daemon="${TRUEPREFIX}/sbin/cherokee"
daemon_flags="-d"
d9 2
a10 1
pexp="${daemon}[^-]"
@


1.4
log
@Where a program with an rc script _requires_ a certain command line
parameter to daemonize, move the parameter from daemon to daemon_flags,
so that the user cannot inadvertently prevent it from daemonizing by
adjusting the flags.

Discussed with ajacoutot and schwarze, this method was suggested
by schwarze@@ as a simpler alternative to my diff. ok aja@@
@
text
@d3 1
a3 1
# $OpenBSD: cherokee.rc,v 1.3 2010/12/27 14:50:24 ajacoutot Exp $
d5 2
a6 1
daemon="${TRUEPREFIX}/sbin/cherokee -d"
@


1.3
log
@Simplify after recent rc.subr change.
The framework is now stable and we will start documenting it (at last).
@
text
@d3 1
a3 1
# $OpenBSD: cherokee.rc,v 1.2 2010/12/24 10:40:06 ajacoutot Exp $
d5 1
a5 2
daemon="${TRUEPREFIX}/sbin/cherokee"
daemon_flags="-d"
@


1.2
log
@Cope with recent rc.subr changes.
@
text
@d3 1
a3 3
# $OpenBSD: cherokee.rc,v 1.1 2010/12/20 07:51:58 ajacoutot Exp $

. /etc/rc.d/rc.subr
d8 1
a8 1
rc_conf
@


1.1
log
@Update to cherokee-1.0.14.

From Fernando Quintero (maintainer).
Several fixes/tweaks as well as the rc script added by myself.
@
text
@d3 1
a3 1
# $OpenBSD$
d9 3
@

