# cvs/ports/www/gitea/pkg/gitea.rc,v
head	1.3;
access;
symbols
	OPENBSD_6_2:1.3.0.2
	OPENBSD_6_2_BASE:1.3
	OPENBSD_6_1:1.1.1.1.0.2
	OPENBSD_6_1_BASE:1.1.1.1
	robert_20170108:1.1.1.1
	robert:1.1.1;
locks; strict;
comment	@ * @;


1.3
date	2017.05.27.19.57.07;	author ajacoutot;	state Exp;
branches;
next	1.2;
commitid	H08QxD8pqoInz6Rz;

1.2
date	2017.05.27.13.53.07;	author ajacoutot;	state Exp;
branches;
next	1.1;
commitid	SDpZdO8hfYIlMQ9j;

1.1
date	2017.01.08.19.20.33;	author robert;	state Exp;
branches
	1.1.1.1;
next	;
commitid	KgfLNVXNsfk70c69;

1.1.1.1
date	2017.01.08.19.20.33;	author robert;	state Exp;
branches;
next	;
commitid	KgfLNVXNsfk70c69;


desc
@@


1.3
log
@Bring rc_bg back.
@
text
@#!/bin/sh
#
# $OpenBSD: gitea.rc,v 1.2 2017/05/27 13:53:07 ajacoutot Exp $

daemon="${TRUEPREFIX}/sbin/gitea"
daemon_user="_gitea"
daemon_flags="web"

. /etc/rc.d/rc.subr

rc_bg=YES
rc_reload=NO

rc_start() {
	${rcexec} "env HOME=${LOCALSTATEDIR}/gitea USER=${daemon_user} \
		GITEA_WORK_DIR=${TRUEPREFIX}/share/gitea \
		GITEA_CUSTOM=${SYSCONFDIR}/gitea \
		${daemon} ${daemon_flags}"
}

rc_cmd $1
@


1.2
log
@Cope with recent rc.subr(8) change (removal of rc_bg etc.)
@
text
@d3 1
a3 1
# $OpenBSD: gitea.rc,v 1.1.1.1 2017/01/08 19:20:33 robert Exp $
d11 1
@


1.1
log
@Initial revision
@
text
@d3 1
a3 1
# $OpenBSD$
a10 1
rc_bg=YES
d17 1
a17 1
		${daemon} ${daemon_flags} ${_bg}"
@


1.1.1.1
log
@import of gitea-1.0.1:

Gitea - Git with a cup of tea is a painless self-hosted Git service comparable 
to Github, Bitbucket or Gitlab.

ok aja@@
@
text
@@
