# cvs/ports/www/gitea/pkg/gitea.rc,v
head     1.1;
branch   1.1.1;
access   ;
symbols  robert_20170108:1.1.1.1 robert:1.1.1;
locks    ; strict;
comment  @ * @;


1.1
date     2017.01.08.19.20.33;  author robert;  state Exp;
branches 1.1.1.1;
next     ;
commitid        KgfLNVXNsfk70c69;

1.1.1.1
date     2017.01.08.19.20.33;  author robert;  state Exp;
branches ;
next     ;
commitid        KgfLNVXNsfk70c69;


desc
@@



1.1
log
@Initial revision
@
text
@#!/bin/sh
#
# $OpenBSD$

daemon="${TRUEPREFIX}/sbin/gitea"
daemon_user="_gitea"
daemon_flags="web"

. /etc/rc.d/rc.subr

rc_bg=YES
rc_reload=NO

rc_start() {
	${rcexec} "env HOME=${LOCALSTATEDIR}/gitea USER=${daemon_user} \
		GITEA_WORK_DIR=${TRUEPREFIX}/share/gitea \
		GITEA_CUSTOM=${SYSCONFDIR}/gitea \
		${daemon} ${daemon_flags} ${_bg}"
}

rc_cmd $1
@


1.1.1.1
log
@import of gitea-1.0.1:

Gitea - Git with a cup of tea is a painless self-hosted Git service comparable 
to Github, Bitbucket or Gitlab.

ok aja@@
@
text
@@
