# cvs/ports/databases/pgcluu/pkg/pgcluu_collectd.rc,v
head	1.2;
access;
symbols
	OPENBSD_6_1:1.2.0.8
	OPENBSD_6_1_BASE:1.2
	OPENBSD_6_0:1.2.0.6
	OPENBSD_6_0_BASE:1.2
	OPENBSD_5_9:1.2.0.2
	OPENBSD_5_9_BASE:1.2
	OPENBSD_5_8:1.2.0.4
	OPENBSD_5_8_BASE:1.2
	jasper_20151706:1.1.1.1
	jasper:1.1.1;
locks; strict;
comment	@ * @;


1.2
date	2015.06.17.12.10.42;	author jasper;	state Exp;
branches;
next	1.1;
commitid	tomvVmH9j734l8Fc;

1.1
date	2015.06.17.10.06.34;	author jasper;	state Exp;
branches
	1.1.1.1;
next	;
commitid	vSMNWp33SOJGEBqr;

1.1.1.1
date	2015.06.17.10.06.34;	author jasper;	state Exp;
branches;
next	;
commitid	vSMNWp33SOJGEBqr;


desc
@@


1.2
log
@start to fill out some sysinfo related fields
@
text
@#!/bin/sh
#
# $OpenBSD: pgcluu_collectd.rc,v 1.1.1.1 2015/06/17 10:06:34 jasper Exp $

daemon="${TRUEPREFIX}/bin/pgcluu_collectd --daemon"
daemon_flags="-i 60 -h localhost -U postgres --disable-sar -f /var/run/pgcluu_collectd.pid /var/db/pgcluu"

. /etc/rc.d/rc.subr

pexp="perl: pgcluu_collectd"
rc_reload=NO

rc_pre() {
	install -d /var/db/pgcluu
	# it doesn't clean up after itself and fails to start if pid file is found.
	rm -f /var/run/pgcluu_collectd.pid
}

rc_cmd $1
@


1.1
log
@Initial revision
@
text
@d3 1
a3 1
# $OpenBSD$
d6 1
a6 1
daemon_flags="-i 60 -h localhost -U postgres --disable-sar --no-sysinfo -f /var/run/pgcluu_collectd.pid /var/db/pgcluu"
@


1.1.1.1
log
@import pgcluu-2.3

PgCluu is a PostgreSQL performances monitoring and auditing tool.  View
reports of all statistics collected from your PostgreSQL databases
cluster. pgCluu will show you the entire informations of the PostgreSQL
Cluster and the system utilization.

ok aja@@
@
text
@@
