# cvs/ports/net/snort/pkg/snort.rc,v
head	1.2;
access;
symbols
	OPENBSD_6_2:1.2.0.6
	OPENBSD_6_2_BASE:1.2
	OPENBSD_6_1:1.2.0.4
	OPENBSD_6_1_BASE:1.2
	OPENBSD_6_0:1.2.0.2
	OPENBSD_6_0_BASE:1.2
	OPENBSD_5_9:1.1.0.12
	OPENBSD_5_9_BASE:1.1
	OPENBSD_5_8:1.1.0.14
	OPENBSD_5_8_BASE:1.1
	OPENBSD_5_7:1.1.0.10
	OPENBSD_5_7_BASE:1.1
	OPENBSD_5_6:1.1.0.8
	OPENBSD_5_6_BASE:1.1
	OPENBSD_5_5:1.1.0.6
	OPENBSD_5_5_BASE:1.1
	OPENBSD_5_4:1.1.0.4
	OPENBSD_5_4_BASE:1.1
	OPENBSD_5_3:1.1.0.2
	OPENBSD_5_3_BASE:1.1;
locks; strict;
comment	@ * @;


1.2
date	2016.05.02.02.29.58;	author lteo;	state Exp;
branches;
next	1.1;
commitid	wdV4pDvleLoy2p6I;

1.1
date	2012.10.11.02.40.48;	author lteo;	state Exp;
branches;
next	;


desc
@@


1.2
log
@Update to Snort 2.9.8.2.

Other minor changes:
- Shortened a line in pkg/README so that portcheck won't complain.
- Replaced "/var" with ${LOCALSTATEDIR} in pkg/snort.rc.
- Regenerated patches.

Tested on amd64 by myself and on i386 by Markus Lude.
Tests were done with DAQ 2.0.6.

ok Markus Lude (maintainer)
@
text
@#!/bin/sh
#
# $OpenBSD: snort.rc,v 1.1 2012/10/11 02:40:48 lteo Exp $

daemon="${TRUEPREFIX}/bin/snort -D"
daemon_flags="-c ${SYSCONFDIR}/snort/snort.conf -u _snort -g _snort -t ${LOCALSTATEDIR}/snort -l ${LOCALSTATEDIR}/snort/log"

. /etc/rc.d/rc.subr

rc_cmd $1
@


1.1
log
@Make various improvements to the Snort port:

* Add an rc.d script.

* In snort.conf, provide the URL to the official Snort rules so that
  users know where to get them.

* In snort.conf, provide the URL to the Emerging Threats rules along
  with a commented include line to allow users to easily load the
  Emerging Threats rules if they wish.

* Revise pkg/README with details on where to obtain Snort rules, the
  differences between the official Snort rules and Emerging Threats
  rules, how to download them, and provide some guidance on setting up
  Snort.

snort.conf and README changes OK Markus Lude (maintainer), sthen@@
rc.d script OK sthen@@
@
text
@d3 1
a3 1
# $OpenBSD$
d6 1
a6 1
daemon_flags="-c ${SYSCONFDIR}/snort/snort.conf -u _snort -g _snort -t /var/snort -l /var/snort/log"
@

