# cvs/ports/www/havp/pkg/havp.rc,v
head	1.2;
access;
symbols
	OPENBSD_6_2:1.2.0.24
	OPENBSD_6_2_BASE:1.2
	OPENBSD_6_1:1.2.0.22
	OPENBSD_6_1_BASE:1.2
	OPENBSD_6_0:1.2.0.20
	OPENBSD_6_0_BASE:1.2
	OPENBSD_5_9:1.2.0.16
	OPENBSD_5_9_BASE:1.2
	OPENBSD_5_8:1.2.0.18
	OPENBSD_5_8_BASE:1.2
	OPENBSD_5_7:1.2.0.14
	OPENBSD_5_7_BASE:1.2
	OPENBSD_5_6:1.2.0.12
	OPENBSD_5_6_BASE:1.2
	OPENBSD_5_5:1.2.0.10
	OPENBSD_5_5_BASE:1.2
	OPENBSD_5_4:1.2.0.8
	OPENBSD_5_4_BASE:1.2
	OPENBSD_5_3:1.2.0.6
	OPENBSD_5_3_BASE:1.2
	OPENBSD_5_2:1.2.0.4
	OPENBSD_5_2_BASE:1.2
	OPENBSD_5_1_BASE:1.2
	OPENBSD_5_1:1.2.0.2
	OPENBSD_5_0:1.1.0.4
	OPENBSD_5_0_BASE:1.1
	OPENBSD_4_9:1.1.0.2
	OPENBSD_4_9_BASE:1.1;
locks; strict;
comment	@ * @;


1.2
date	2011.11.30.11.00.59;	author robert;	state Exp;
branches;
next	1.1;

1.1
date	2010.12.30.10.21.25;	author giovanni;	state Exp;
branches;
next	;


desc
@@


1.2
log
@create the working directories with the correct group
so that others can access it
discussed with ajacoutot@@
@
text
@#!/bin/sh
#
# $OpenBSD: havp.rc,v 1.1 2010/12/30 10:21:25 giovanni Exp $

daemon="${TRUEPREFIX}/sbin/havp"

. /etc/rc.d/rc.subr

rc_pre() {
	install -d -o _havp -g _havp /var/run/havp/ /var/tmp/havp/
}

rc_cmd $1
@


1.1
log
@
Update to 0.92a and add an rcscript
Hints and ok ajacoutot@@, ok benoit@@
@
text
@d3 1
a3 1
# $OpenBSD$
d10 1
a10 1
	install -d -o _havp /var/run/havp/ /var/tmp/havp/
@

