# cvs/ports/security/cfs/pkg/Attic/PLIST,v
head	1.8;
access;
symbols
	OPENBSD_5_8:1.7.0.46
	OPENBSD_5_8_BASE:1.7
	OPENBSD_5_7:1.7.0.42
	OPENBSD_5_7_BASE:1.7
	OPENBSD_5_6:1.7.0.40
	OPENBSD_5_6_BASE:1.7
	OPENBSD_5_5:1.7.0.38
	OPENBSD_5_5_BASE:1.7
	OPENBSD_5_4:1.7.0.36
	OPENBSD_5_4_BASE:1.7
	OPENBSD_5_3:1.7.0.34
	OPENBSD_5_3_BASE:1.7
	OPENBSD_5_2:1.7.0.32
	OPENBSD_5_2_BASE:1.7
	OPENBSD_5_1_BASE:1.7
	OPENBSD_5_1:1.7.0.30
	OPENBSD_5_0:1.7.0.28
	OPENBSD_5_0_BASE:1.7
	OPENBSD_4_9:1.7.0.26
	OPENBSD_4_9_BASE:1.7
	OPENBSD_4_8:1.7.0.24
	OPENBSD_4_8_BASE:1.7
	OPENBSD_4_7:1.7.0.22
	OPENBSD_4_7_BASE:1.7
	OPENBSD_4_6:1.7.0.20
	OPENBSD_4_6_BASE:1.7
	OPENBSD_4_5:1.7.0.18
	OPENBSD_4_5_BASE:1.7
	OPENBSD_4_4:1.7.0.16
	OPENBSD_4_4_BASE:1.7
	OPENBSD_4_3:1.7.0.14
	OPENBSD_4_3_BASE:1.7
	OPENBSD_4_2:1.7.0.12
	OPENBSD_4_2_BASE:1.7
	OPENBSD_4_1:1.7.0.10
	OPENBSD_4_1_BASE:1.7
	OPENBSD_4_0:1.7.0.8
	OPENBSD_4_0_BASE:1.7
	OPENBSD_3_9:1.7.0.6
	OPENBSD_3_9_BASE:1.7
	OPENBSD_3_8:1.7.0.4
	OPENBSD_3_8_BASE:1.7
	OPENBSD_3_7:1.7.0.2
	OPENBSD_3_7_BASE:1.7
	OPENBSD_3_6:1.5.0.2
	OPENBSD_3_6_BASE:1.5
	OPENBSD_3_5:1.4.0.10
	OPENBSD_3_5_BASE:1.4
	OPENBSD_3_4:1.4.0.8
	OPENBSD_3_4_BASE:1.4
	OPENBSD_3_3:1.4.0.6
	OPENBSD_3_3_BASE:1.4
	OPENBSD_3_2:1.4.0.4
	OPENBSD_3_2_BASE:1.4
	OPENBSD_3_1:1.4.0.2
	OPENBSD_3_1_BASE:1.4
	OPENBSD_3_0:1.3.0.8
	OPENBSD_3_0_BASE:1.3
	OPENBSD_2_9_TRACKING_SWITCH:1.3
	OPENBSD_2_9:1.3.0.6
	OPENBSD_2_9_BASE:1.3
	OPENBSD_2_8:1.3.0.4
	OPENBSD_2_8_BASE:1.3
	OPENBSD_2_7:1.3.0.2
	OPENBSD_2_7_BASE:1.3
	OPENBSD_2_6:1.2.0.4
	OPENBSD_2_6_BASE:1.2
	OPENBSD_2_5:1.2.0.2
	OPENBSD_2_5_BASE:1.2
	OPENBSD_2_4:1.1.0.4
	OPENBSD_2_4_BASE:1.1
	OPENBSD_2_3:1.1.0.2
	OPENBSD_2_3_BASE:1.1;
locks; strict;
comment	@# @;


1.8
date	2015.08.25.11.16.28;	author jca;	state dead;
branches;
next	1.7;
commitid	ZYz3xUSOJVyi0xdm;

1.7
date	2004.12.03.16.08.48;	author alek;	state Exp;
branches;
next	1.6;

1.6
date	2004.09.15.18.35.59;	author espie;	state Exp;
branches;
next	1.5;

1.5
date	2004.06.15.18.32.22;	author sturm;	state Exp;
branches;
next	1.4;

1.4
date	2002.04.07.01.54.52;	author naddy;	state Exp;
branches;
next	1.3;

1.3
date	99.12.09.04.46.18;	author brad;	state Exp;
branches;
next	1.2;

1.2
date	99.02.28.08.35.43;	author fgsch;	state Exp;
branches;
next	1.1;

1.1
date	97.11.27.22.22.53;	author angelos;	state Exp;
branches;
next	;


desc
@@


1.8
log
@Remove security/cfs.  Unmaintained code, broken on 64 bits archs.

softraid CRYPTO, vnconfig and security/encfs provide much saner
alternatives.

ok tedu@@ steven@@ ajacoutot@@ sthen@@
@
text
@@@comment $OpenBSD: PLIST,v 1.7 2004/12/03 16:08:48 alek Exp $
bin/cattach
bin/cdetach
bin/cmkdir
bin/cmkkey
bin/cpasswd
@@man man/man1/cattach.1
@@man man/man1/cdetach.1
@@man man/man1/cmkdir.1
@@man man/man1/cmkkey.1
@@man man/man1/cpasswd.1
@@man man/man8/cfsd.8
sbin/cfsd
share/doc/cfs/
share/doc/cfs/README.install
share/doc/cfs/notes.ms
@@sample /crypt/
@@mode 0
@@sample /null/
@@extraunexec umount /crypt > /dev/null 2>&1 || true
@@extraunexec rm -r /crypt
@@extraunexec rm -r /null
@


1.7
log
@- Replace INSTALL with @@sample
- Better @@extraunexec
- Replace DEINSTALL with UNMESSAGE
- Bump PKGNAME
- Add WANTLIB marker

ok naddy@@
@
text
@d1 1
a1 1
@@comment $OpenBSD: PLIST,v 1.6 2004/09/15 18:35:59 espie Exp $
@


1.6
log
@new plists.
gnutls depends on lzo.
@
text
@d1 1
a1 1
@@comment $OpenBSD: PLIST,v 1.5 2004/06/15 18:32:22 sturm Exp $
d17 4
a20 1
@@extraunexec umount /crypt > /dev/null 2>&1
@


1.5
log
@update to cfs 1.4.1

CFS was actually broken since at least 3.2 or so, but should work now
again on 32-bit archs. 64-bit archs are still broken and won't be
fixed anytime soon.

I also set PERMIT_*=Yes, as the "no fee" part of the license is clearly
due to bad english.

xsa@@ helped testing
@
text
@d1 1
a1 1
@@comment $OpenBSD: PLIST,v 1.4 2002/04/07 01:54:52 naddy Exp $
d7 6
a12 6
man/man1/cattach.1
man/man1/cdetach.1
man/man1/cmkdir.1
man/man1/cmkkey.1
man/man1/cpasswd.1
man/man8/cfsd.8
d14 1
a16 1
@@dirrm share/doc/cfs
@


1.4
log
@Add missing RCS Id.
@
text
@d1 1
a1 1
@@comment $OpenBSD$
d5 1
d10 1
a11 1
man/man8/ccat.8
a12 2
man/man8/cname.8
sbin/ccat
a13 1
sbin/cname
d17 3
@


1.3
log
@- fix install target
- add missing @@dirrm to PLIST
@
text
@d1 1
@


1.2
log
@Don't install cfs.sh, it is started in /etc/rc.local.
@
text
@d17 1
@


1.1
log
@CFS 1.3.3 + blowfish (from Niels Provos) port.
It's fun being in Canada.
@
text
@a4 1
etc/rc.d/cfsd.sh
@

