# cvs/ports/net/smokeping/pkg/smokeping_fcgi.rc,v
head	1.1;
access;
symbols
	OPENBSD_6_1:1.1.0.10
	OPENBSD_6_1_BASE:1.1
	OPENBSD_6_0:1.1.0.8
	OPENBSD_6_0_BASE:1.1
	OPENBSD_5_9:1.1.0.4
	OPENBSD_5_9_BASE:1.1
	OPENBSD_5_8:1.1.0.6
	OPENBSD_5_8_BASE:1.1
	OPENBSD_5_7:1.1.0.2
	OPENBSD_5_7_BASE:1.1;
locks; strict;
comment	@ * @;


1.1
date	2014.11.05.15.42.09;	author sthen;	state Exp;
branches;
next	;
commitid	r3Gb5UT6Agg7kLiO;


desc
@@


1.1
log
@update to smokeping 2.6.11
- add a new rc script to run as FastCGI
@
text
@#!/bin/sh
#
# $OpenBSD: smokeping.rc,v 1.1 2013/01/05 22:52:07 sthen Exp $

daemon="${TRUEPREFIX}/bin/smokeping_cgi"
daemon_flags="${SYSCONFDIR}/config"

. /etc/rc.d/rc.subr

pexp="/usr/bin/perl ${daemon}${daemon_flags:+ ${daemon_flags}}"
rc_reload=NO

rc_start() {
	${LOCALBASE}/bin/spawn-fcgi -s /var/www/run/smokeping.sock \
	    -u _smokeping -U www -M 0660 ${daemon} ${daemon_flags}
}

rc_cmd $1
@
