# cvs/ports/www/tt-rss/pkg/Attic/tt_rss_update.rc,v
head	1.6;
access;
symbols
	OPENBSD_5_7:1.5.0.8
	OPENBSD_5_7_BASE:1.5
	OPENBSD_5_6:1.5.0.6
	OPENBSD_5_6_BASE:1.5
	OPENBSD_5_5:1.5.0.4
	OPENBSD_5_5_BASE:1.5
	OPENBSD_5_4:1.5.0.2
	OPENBSD_5_4_BASE:1.5
	OPENBSD_5_3:1.4.0.2
	OPENBSD_5_3_BASE:1.4
	OPENBSD_5_2:1.3.0.2
	OPENBSD_5_2_BASE:1.3;
locks; strict;
comment	@ * @;


1.6
date	2015.06.26.14.00.58;	author ajacoutot;	state dead;
branches;
next	1.5;
commitid	EGyW970RANpepzVp;

1.5
date	2013.04.14.20.28.09;	author rpe;	state Exp;
branches;
next	1.4;

1.4
date	2012.08.25.11.41.46;	author ajacoutot;	state Exp;
branches;
next	1.3;

1.3
date	2012.05.03.18.45.39;	author landry;	state Exp;
branches;
next	1.2;

1.2
date	2012.03.15.11.28.19;	author ajacoutot;	state Exp;
branches;
next	1.1;

1.1
date	2012.03.15.09.50.23;	author landry;	state Exp;
branches;
next	;


desc
@@


1.6
log
@Remove. Upstream will not make releases anymore and advise to use git
master instead. It's pretty easy to do a checkout and run it instead of
using a package.

ok jasper@@
@
text
@#!/bin/sh
#
# $OpenBSD: tt_rss_update.rc,v 1.5 2013/04/14 20:28:09 rpe Exp $

daemon="${INSTDIR}/update.php --daemon"
daemon_user="www"

. /etc/rc.d/rc.subr

pexp="${MODPHP_BIN} ${daemon}${daemon_flags:+ ${daemon_flags}}"

rc_bg=YES
rc_reload=NO

rc_post() {
	rm -f ${INSTDIR}/lock/update_daemon.lock
}

rc_cmd $1
@


1.5
log
@update tt-rss to 1.7.8

with feedback and ok landry@@ sthen@@
@
text
@d3 1
a3 1
# $OpenBSD: tt_rss_update.rc,v 1.4 2012/08/25 11:41:46 ajacoutot Exp $
@


1.4
log
@Make sure rc_post does not fail by using rm -f ; lockfiles aren't
supposed to stay around.

from rte
@
text
@d3 1
a3 1
# $OpenBSD: tt_rss_update.rc,v 1.3 2012/05/03 18:45:39 landry Exp $
d5 1
a5 1
daemon="${INSTDIR}/update.php -daemon"
@


1.3
log
@Fix rc script so that rc_post does what it's supposed to do.
@
text
@d3 1
a3 1
# $OpenBSD: tt_rss_update.rc,v 1.2 2012/03/15 11:28:19 ajacoutot Exp $
d16 1
a16 1
	rm ${INSTDIR}/lock/update_daemon.lock
@


1.2
log
@* don't append a trailing '/' to INSTDIR
* fix crontab entry
* sort rc script

ok landry@@
@
text
@d3 1
a3 1
# $OpenBSD: tt_rss_update.rc,v 1.1 2012/03/15 09:50:23 landry Exp $
d16 1
a16 1
	rm ${INSTDIR}lock/update_daemon.lock
@


1.1
log
@Update to tt-rss 1.5.10

- remove dep on gtar
- patch config.php-dist & update.php to use MODPHP_BIN
- fix perms in PLIST, use @@owner www for dirs where tt-rss writes
- switch MESSAGE to README, and improve the latter
- provide an rc script for the daemon updating the feeds

ok ajacoutot@@
@
text
@d3 1
a3 1
# $OpenBSD$
d5 1
a5 1
daemon="${INSTDIR}update.php -daemon"
d12 1
a13 1
rc_bg=YES
@

