# cvs/ports/net/bitlbee/pkg/bitlbee.rc,v
head	1.1;
access;
symbols
	OPENBSD_6_2:1.1.0.2
	OPENBSD_6_2_BASE:1.1;
locks; strict;
comment	@ * @;


1.1
date	2017.05.29.12.11.07;	author sthen;	state Exp;
branches;
next	;
commitid	7gcQJTkWdnLWtzke;


desc
@@


1.1
log
@add an rc script for bitlbee, ok Tom Doherty (maintainer), ajacoutot@@
@
text
@#!/bin/sh
#
# $OpenBSD$

daemon="${TRUEPREFIX}/sbin/bitlbee -F"
daemon_user="_bitlbee"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
@
