# cvs/ports/security/ssh/pkg/Attic/PLIST,v
head	1.11;
access;
symbols
	OPENBSD_2_7:1.10.0.6
	OPENBSD_2_7_BASE:1.10
	OPENBSD_2_6:1.10.0.4
	OPENBSD_2_6_BASE:1.10
	OPENBSD_2_5:1.10.0.2
	OPENBSD_2_5_BASE:1.10
	OPENBSD_2_4:1.9.0.2
	OPENBSD_2_4_BASE:1.9
	OPENBSD_2_3:1.6.0.2
	OPENBSD_2_3_BASE:1.6;
locks; strict;
comment	@# @;


1.11
date	2000.08.27.20.45.16;	author turan;	state dead;
branches;
next	1.10;

1.10
date	99.03.23.17.48.05;	author espie;	state Exp;
branches;
next	1.9;

1.9
date	98.07.13.01.17.02;	author marc;	state Exp;
branches;
next	1.8;

1.8
date	98.07.07.04.00.00;	author marc;	state Exp;
branches;
next	1.7;

1.7
date	98.06.02.18.10.12;	author marc;	state Exp;
branches;
next	1.6;

1.6
date	98.04.05.07.29.07;	author marc;	state Exp;
branches;
next	1.5;

1.5
date	98.02.23.19.50.23;	author don;	state Exp;
branches;
next	1.4;

1.4
date	98.02.13.14.46.00;	author niklas;	state Exp;
branches;
next	1.3;

1.3
date	98.01.23.05.12.51;	author todd;	state Exp;
branches;
next	1.2;

1.2
date	97.11.25.22.59.46;	author angelos;	state Exp;
branches;
next	1.1;

1.1
date	97.11.17.05.58.47;	author todd;	state Exp;
branches;
next	;


desc
@@


1.11
log
@remove ssh port.  discussed with theo.
@
text
@bin/ssh1
bin/ssh
bin/slogin
bin/ssh-keygen1
bin/ssh-keygen
bin/ssh-agent1
bin/ssh-agent
bin/ssh-add1
bin/ssh-add
bin/scp1
bin/scp
bin/ssh-askpass1
bin/ssh-askpass
bin/make-ssh-known-hosts1
bin/make-ssh-known-hosts
lib/ssh/server_config.sample
@@exec if [ ! -f /etc/sshd_config ]; then cp %D/%F /etc/sshd_config; fi
lib/ssh/host_config.sample
@@exec if [ ! -f /etc/ssh_config ]; then cp %D/%F /etc/ssh_config; fi
man/man1/ssh-keygen1.1
man/man1/ssh-keygen.1
man/man1/ssh-agent1.1
man/man1/ssh-agent.1
man/man1/ssh-add1.1
man/man1/ssh-add.1
man/man1/scp1.1
man/man1/scp.1
man/man1/slogin1.1
man/man1/slogin.1
man/man1/ssh1.1
man/man1/ssh.1
man/man1/make-ssh-known-hosts1.1
man/man1/make-ssh-known-hosts.1
man/man8/sshd1.8
man/man8/sshd.8
sbin/sshd1
sbin/sshd
@@dirrm lib/ssh
@@exec if [ ! -f /etc/ssh_host_key ]; then echo "Generating a secret host key.." ; %D/bin/ssh-keygen -N "" -f /etc/ssh_host_key; fi
@@unexec echo "***"
@@unexec echo "*** ssh configuration files stored in /etc have not been removed."
@@unexec echo "*** you may wish to remove them manually."
@@unexec echo "***"
@


1.10
log
@Fix up PLISTs.

At one point bsd.port.mk put out erroneous PLIST-auto, with @@cwd, @@pkgname,
and @@pkgdep lines.

These are actually added by make package, and shouldn't usually be put
inside the PLIST, as this is extra information to update.

I left a few @@cwd in, mainly for the perl modules, which want to live under
/usr to work.
@
text
@@


1.9
log
@yet another version of ssh -- 1.2.26
@
text
@a0 2
@@cwd /usr/local
@@name ssh-1.2.26
@


1.8
log
@add patch from http://www.core-sdi.com/ssh/
@
text
@a0 2
@@comment PACKAGE(arch=sparc, opsys=OpenBSD, vers=2.3)
@@name ssh-1.2.25
d2 1
@


1.7
log
@update to ssh 1.2.23
@
text
@d1 3
d5 2
a6 4
@@exec ln -fs %f %B/ssh
@@unexec rm -f %B/ssh
@@exec ln -fs %B/ssh %B/slogin
@@unexec rm -f %B/slogin
d8 1
a8 2
@@exec ln -fs %f %B/ssh-keygen
@@unexec rm -f %B/ssh-keygen
d10 1
a10 2
@@exec ln -fs %f %B/ssh-agent
@@unexec rm -f %B/ssh-agent
d12 1
a12 2
@@exec ln -fs %f %B/ssh-add
@@unexec rm -f %B/ssh-add
d14 1
a14 2
@@exec ln -fs %f %B/scp
@@unexec rm -f %B/scp
d16 1
a16 2
@@exec ln -fs %f %B/ssh-askpass
@@unexec rm -f %B/ssh-askpass
d18 5
a22 11
@@exec ln -fs %f %B/make-ssh-known-hosts
@@unexec rm -f %B/make-ssh-known-hosts
sbin/sshd1
@@exec ln -fs %f %B/sshd
@@unexec rm -f %B/sshd
man/man1/make-ssh-known-hosts1.1
@@exec ln -fs %f %B/make-ssh-known-hosts.1
@@unexec rm -f %B/make-ssh-known-hosts.1
man/man1/scp1.1
@@exec ln -fs %f %B/scp.1
@@unexec rm -f %B/scp.1
d24 1
a24 2
@@exec ln -fs %f %B/ssh-keygen.1
@@unexec rm -f %B/ssh-keygen.1
d26 1
a26 2
@@exec ln -fs %f %B/ssh-agent.1
@@unexec rm -f %B/ssh-agent.1
d28 5
a32 2
@@exec ln -fs %f %B/ssh-add.1
@@unexec rm -f %B/ssh-add.1
d34 3
a36 6
@@exec ln -fs %f %B/ssh.1
@@unexec rm -f %B/ssh.1
@@exec ln -fs %f %B/slogin1.1
@@unexec rm -f %B/slogin1.1
@@exec ln -fs %B/ssh.1 %B/slogin.1
@@unexec rm -f %B/slogin.1
d38 4
a41 6
@@exec ln -fs %f %B/sshd.8
@@unexec rm -f %B/sshd.8
lib/ssh/server_config.sample
@@exec if [ ! -f /etc/sshd_config ]; then cp %D/%F /etc/sshd_config; fi
lib/ssh/host_config.sample
@@exec if [ ! -f /etc/ssh_config ]; then cp %D/%F /etc/ssh_config; fi
a42 1
@@dirrm lib/ssh
@


1.6
log
@Save sample config files in a place where they can be referenced if needed
when installing from a package.  Sample files will NOT replace any existing
configuration files.
@
text
@d1 4
a4 3
bin/scp
bin/ssh
@@exec ln -fs %f %B/slogin
d6 42
a47 12
bin/ssh-add
bin/ssh-agent
bin/ssh-askpass
bin/ssh-keygen
bin/make-ssh-known-hosts
man/man1/make-ssh-known-hosts.1
man/man1/scp.1
man/man1/ssh-add.1
man/man1/ssh-agent.1
man/man1/ssh-keygen.1
man/man1/ssh.1
@@exec ln -fs %f %B/slogin.1
d49 3
a51 2
man/man8/sshd.8
sbin/sshd
d58 4
a61 2
@@unexec echo "ssh configuration files stored in /etc have not been removed."
@@unexec echo "you may wish to remove them manually."
@


1.5
log
@include tcp wrappers automatically; misc cleanup
@
text
@d20 4
d25 1
@


1.4
log
@We do not have host-specific things in /usr/local/etc
@
text
@a0 1
etc/rc.d/sshd.sh
d21 2
@


1.3
log
@we use /etc not /usr/local/etc for ssh config files
@
text
@a10 2
etc/ssh_config
etc/sshd_config
@


1.2
log
@Fix PLIST for man pages.
@
text
@d23 1
a23 1
@@exec if [ ! -f %D/etc/ssh_host_key ]; then echo "Generating a secret host key.." ; %D/bin/ssh-keygen -N "" -f %D/etc/ssh_host_key; fi
@


1.1
log
@+ ssh (1.2.21)
@
text
@d13 9
a21 9
man/man1/make-ssh-known-hosts.1.gz
man/man1/scp.1.gz
man/man1/ssh-add.1.gz
man/man1/ssh-agent.1.gz
man/man1/ssh-keygen.1.gz
man/man1/ssh.1.gz
@@exec ln -fs %f %B/slogin.1.gz
@@unexec rm -f %B/slogin.1.gz
man/man8/sshd.8.gz
@
