# cvs/ports/net/lldpd/pkg/lldpd.rc,v
head	1.4;
access;
symbols
	OPENBSD_6_2:1.4.0.6
	OPENBSD_6_2_BASE:1.4
	OPENBSD_6_1:1.4.0.4
	OPENBSD_6_1_BASE:1.4
	OPENBSD_6_0:1.4.0.2
	OPENBSD_6_0_BASE:1.4
	OPENBSD_5_9:1.2.0.10
	OPENBSD_5_9_BASE:1.2
	OPENBSD_5_8:1.2.0.12
	OPENBSD_5_8_BASE:1.2
	OPENBSD_5_7:1.2.0.8
	OPENBSD_5_7_BASE:1.2
	OPENBSD_5_6:1.2.0.6
	OPENBSD_5_6_BASE:1.2
	OPENBSD_5_5:1.2.0.4
	OPENBSD_5_5_BASE:1.2
	OPENBSD_5_4:1.2.0.2
	OPENBSD_5_4_BASE:1.2
	sthen_20130712:1.1.1.1
	sthen:1.1.1;
locks; strict;
comment	@ * @;


1.4
date	2016.06.28.17.03.27;	author ajacoutot;	state Exp;
branches;
next	1.3;
commitid	2Ty2T8CPfuAPg7SR;

1.3
date	2016.06.28.15.50.27;	author sthen;	state Exp;
branches;
next	1.2;
commitid	PZqMZh0aHsrVpRYL;

1.2
date	2013.07.15.21.35.46;	author sthen;	state Exp;
branches;
next	1.1;

1.1
date	2013.07.12.10.26.44;	author sthen;	state Exp;
branches
	1.1.1.1;
next	;

1.1.1.1
date	2013.07.12.10.26.44;	author sthen;	state Exp;
branches;
next	;


desc
@@


1.4
log
@rc_ overrides are put after sourcing /etc/rc.d/rc.subr.

ok sthen@@ (maintainer)
@
text
@#!/bin/sh
#
# $OpenBSD: lldpd.rc,v 1.3 2016/06/28 15:50:27 sthen Exp $

daemon="${TRUEPREFIX}/sbin/lldpd"

. /etc/rc.d/rc.subr

pexp="lldpd: monitor."
rc_reload=NO

rc_pre() {
	install -d -o _lldpd /var/run/lldpd
}

rc_cmd $1
@


1.3
log
@process titles now have an extra "." at the end; fix lldpd.rc
@
text
@d3 1
a3 1
# $OpenBSD: lldpd.rc,v 1.2 2013/07/15 21:35:46 sthen Exp $
a5 1
rc_reload=NO
d10 1
@


1.2
log
@fix rc.d stop/restart (broken pexp), reported by jasper@@
add sigaction diff from upstream
@
text
@d3 1
a3 1
# $OpenBSD: lldpd.rc,v 1.1.1.1 2013/07/12 10:26:44 sthen Exp $
d10 1
a10 1
pexp="lldpd: monitor"
@


1.1
log
@Initial revision
@
text
@d3 1
a3 1
# $OpenBSD$
d9 2
@


1.1.1.1
log
@import ports/net/lldpd, an alternative to ladvd. ok jasper@@

LLDP (Link Layer Discovery Protocol) is an industry standard protocol
designed to supplant proprietary Link-Layer protocols such as
Extreme's EDP (Extreme Discovery Protocol) and CDP (Cisco Discovery
Protocol). The goal of LLDP is to provide an inter-vendor compatible
mechanism to deliver Link-Layer notifications to adjacent network
devices.

lldpd implements both reception and sending. It also implements an
SNMP subagent for net-snmp to get local and remote LLDP
information. The LLDP MIB is partially implemented but the most useful
tables are here. lldpd also partially implements LLDP-MED.

lldpd supports bridge, vlan and bonding. bonding need to be done on
real physical devices, not on bridges, vlans, etc. However, vlans can
be mapped on the bonding device. You can bridge vlan but not add vlans
on bridges. More complex setups may give false results.
@
text
@@
