# cvs/ports/mail/milter-spamd/pkg/milter_spamd.rc,v
head	1.6;
access;
symbols
	OPENBSD_6_2:1.6.0.10
	OPENBSD_6_2_BASE:1.6
	OPENBSD_6_1:1.6.0.8
	OPENBSD_6_1_BASE:1.6
	OPENBSD_6_0:1.6.0.6
	OPENBSD_6_0_BASE:1.6
	OPENBSD_5_9:1.6.0.2
	OPENBSD_5_9_BASE:1.6
	OPENBSD_5_8:1.6.0.4
	OPENBSD_5_8_BASE:1.6
	OPENBSD_5_7:1.5.0.10
	OPENBSD_5_7_BASE:1.5
	OPENBSD_5_6:1.5.0.8
	OPENBSD_5_6_BASE:1.5
	OPENBSD_5_5:1.5.0.6
	OPENBSD_5_5_BASE:1.5
	OPENBSD_5_4:1.5.0.4
	OPENBSD_5_4_BASE:1.5
	OPENBSD_5_3:1.5.0.2
	OPENBSD_5_3_BASE:1.5
	OPENBSD_5_2:1.4.0.8
	OPENBSD_5_2_BASE:1.4
	OPENBSD_5_1_BASE:1.4
	OPENBSD_5_1:1.4.0.6
	OPENBSD_5_0:1.4.0.4
	OPENBSD_5_0_BASE:1.4
	OPENBSD_4_9:1.4.0.2
	OPENBSD_4_9_BASE:1.4;
locks; strict;
comment	@ * @;


1.6
date	2015.07.15.19.47.35;	author ajacoutot;	state Exp;
branches;
next	1.5;
commitid	ZzLBxIjV0SDzTERN;

1.5
date	2012.12.05.09.38.41;	author ajacoutot;	state Exp;
branches;
next	1.4;

1.4
date	2010.12.27.14.50.23;	author ajacoutot;	state Exp;
branches;
next	1.3;

1.3
date	2010.12.24.10.40.06;	author ajacoutot;	state Exp;
branches;
next	1.2;

1.2
date	2010.12.11.12.57.37;	author ajacoutot;	state Exp;
branches;
next	1.1;

1.1
date	2010.12.09.01.47.31;	author okan;	state Exp;
branches;
next	;


desc
@@


1.6
log
@No need for a specific pexp.
@
text
@#!/bin/sh
#
# $OpenBSD: milter_spamd.rc,v 1.5 2012/12/05 09:38:41 ajacoutot Exp $

daemon="${TRUEPREFIX}/libexec/milter-spamd"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
@


1.5
log
@Make the pexp match on "${daemon}"; when using the `-i' option (for regex),
the process name would not match otherwise.

ok sthen@@
@
text
@d3 1
a3 1
# $OpenBSD: milter_spamd.rc,v 1.4 2010/12/27 14:50:23 ajacoutot Exp $
a8 1
pexp="${daemon}"
@


1.4
log
@Simplify after recent rc.subr change.
The framework is now stable and we will start documenting it (at last).
@
text
@d3 1
a3 1
# $OpenBSD: milter_spamd.rc,v 1.3 2010/12/24 10:40:06 ajacoutot Exp $
d9 1
@


1.3
log
@Cope with recent rc.subr changes.
@
text
@d3 1
a3 3
# $OpenBSD: milter_spamd.rc,v 1.2 2010/12/11 12:57:37 ajacoutot Exp $

. /etc/rc.d/rc.subr
d7 1
a7 1
rc_conf
@


1.2
log
@Move these to rc_reload=NO.
@
text
@d3 1
a3 1
# $OpenBSD: milter_spamd.rc,v 1.1 2010/12/09 01:47:31 okan Exp $
d8 3
@


1.1
log
@- add rc script
- no need for groff
- remove MESSAGE since the remaining relevant snippet is in the manpage
(maintainer timeout)

feedback and ok ajacoutot@@
@
text
@d3 1
a3 1
# $OpenBSD$
d8 1
a8 4

rc_reload() {
	rc_err "$0: reload is not supported"
}
@

