# cvs/ports/net/munin/pkg/Attic/munin-node.rc,v
head	1.2;
access;
symbols;
locks; strict;
comment	@ * @;


1.2
date	2010.10.28.20.59.05;	author jasper;	state dead;
branches;
next	1.1;

1.1
date	2010.10.28.11.01.47;	author jasper;	state Exp;
branches;
next	;


desc
@@


1.2
log
@After discussion with robert and ingo, we are renaming the rc scripts from
dash to underscore so that we are consitent with the foo_bar_flags
variable.
@
text
@#!/bin/sh

. /etc/rc.d/rc.subr

daemon="${TRUEPREFIX}/sbin/munin-node"
pexp="perl: ${daemon}"

rc_pre() {
        install -d -o _munin /var/run/munin
}

rc_cmd $1
@


1.1
log
@- add the munin-node rc script

ok sthen@@ (MAINTAINER) robert@@
@
text
@@

