# cvs/ports/multimedia/motion/pkg/motion.rc,v
head	1.2;
access;
symbols
	OPENBSD_6_2:1.2.0.12
	OPENBSD_6_2_BASE:1.2
	OPENBSD_6_1:1.2.0.10
	OPENBSD_6_1_BASE:1.2
	OPENBSD_6_0:1.2.0.8
	OPENBSD_6_0_BASE:1.2
	OPENBSD_5_9:1.2.0.4
	OPENBSD_5_9_BASE:1.2
	OPENBSD_5_8:1.2.0.6
	OPENBSD_5_8_BASE:1.2
	OPENBSD_5_7:1.2.0.2
	OPENBSD_5_7_BASE:1.2
	sthen_20141010:1.1.1.1
	sthen:1.1.1;
locks; strict;
comment	@ * @;


1.2
date	2014.10.11.11.04.20;	author sthen;	state Exp;
branches;
next	1.1;
commitid	OIqRYyaexut4aC1U;

1.1
date	2014.10.10.13.32.05;	author sthen;	state Exp;
branches
	1.1.1.1;
next	;
commitid	MGj66DDvNXJgYAYG;

1.1.1.1
date	2014.10.10.13.32.05;	author sthen;	state Exp;
branches;
next	;
commitid	MGj66DDvNXJgYAYG;


desc
@@


1.2
log
@Patch motion to support -b (background / daemon mode) overriding the
"daemon" config value. Use it in motion.rc to ensure safe startup even
if config has "daemon off".
@
text
@#!/bin/sh
#
# $OpenBSD: motion.rc,v 1.1.1.1 2014/10/10 13:32:05 sthen Exp $

daemon="${TRUEPREFIX}/bin/motion -b"
daemon_user="_motion"

. /etc/rc.d/rc.subr

rc_pre() {
	install -d -o _motion /var/run/motion
}

rc_cmd $1
@


1.1
log
@Initial revision
@
text
@d3 1
a3 1
# $OpenBSD$
d5 1
a5 1
daemon="${TRUEPREFIX}/bin/motion"
@


1.1.1.1
log
@import ports/multimedia/motion, ok/tweak ajacoutot@@

Motion is a C program that monitors the video signal from one or
more V4L2/video(4) cameras or network-based webcams, and is able to
detect if a significant part of the picture has changed. Or in
other words, it can detect motion.
@
text
@@
