# cvs/ports/net/nagios/nrpe/pkg/nrpe.rc,v
head	1.3;
access;
symbols
	OPENBSD_6_0:1.3.0.16
	OPENBSD_6_0_BASE:1.3
	OPENBSD_5_9:1.3.0.12
	OPENBSD_5_9_BASE:1.3
	OPENBSD_5_8:1.3.0.14
	OPENBSD_5_8_BASE:1.3
	OPENBSD_5_7:1.3.0.10
	OPENBSD_5_7_BASE:1.3
	OPENBSD_5_6:1.3.0.8
	OPENBSD_5_6_BASE:1.3
	OPENBSD_5_5:1.3.0.6
	OPENBSD_5_5_BASE:1.3
	OPENBSD_5_4:1.3.0.4
	OPENBSD_5_4_BASE:1.3
	OPENBSD_5_3:1.3.0.2
	OPENBSD_5_3_BASE:1.3
	OPENBSD_5_2:1.2.0.4
	OPENBSD_5_2_BASE:1.2
	OPENBSD_5_1_BASE:1.2
	OPENBSD_5_1:1.2.0.2
	OPENBSD_5_0:1.1.0.4
	OPENBSD_5_0_BASE:1.1
	OPENBSD_4_9:1.1.0.2
	OPENBSD_4_9_BASE:1.1;
locks; strict;
comment	@ * @;


1.3
date	2012.09.05.19.55.57;	author sthen;	state Exp;
branches;
next	1.2;

1.2
date	2011.12.17.22.01.13;	author sthen;	state Exp;
branches;
next	1.1;

1.1
date	2011.01.20.14.04.42;	author sthen;	state Exp;
branches;
next	;


desc
@@


1.3
log
@use a subdir for the pid file to avoid log spam; pointed out by Julien Dhaille
@
text
@#!/bin/sh
#
# $OpenBSD: nrpe.rc,v 1.2 2011/12/17 22:01:13 sthen Exp $

daemon="${TRUEPREFIX}/sbin/nrpe -d"

. /etc/rc.d/rc.subr

rc_pre() {
	install -o _nrpe -d /var/run/nrpe
}

rc_cmd $1
@


1.2
log
@Where a program with an rc script _requires_ a certain command line
parameter to daemonize, move the parameter from daemon to daemon_flags,
so that the user cannot inadvertently prevent it from daemonizing by
adjusting the flags.

Discussed with ajacoutot and schwarze, this method was suggested
by schwarze@@ as a simpler alternative to my diff. ok aja@@
@
text
@d3 1
a3 1
# $OpenBSD: nrpe.rc,v 1.1 2011/01/20 14:04:42 sthen Exp $
d8 4
@


1.1
log
@add an rc script for nrpe; ok ajacoutot@@
@
text
@d3 1
a3 1
# $OpenBSD$
d5 1
a5 2
daemon="${TRUEPREFIX}/sbin/nrpe"
daemon_flags="-d"
@

