# cvs/ports/www/squid/snapshot/pkg/Attic/squid.rc,v
head	1.4;
access;
symbols
	OPENBSD_5_8:1.3.0.6
	OPENBSD_5_8_BASE:1.3
	OPENBSD_5_7:1.3.0.2
	OPENBSD_5_7_BASE:1.3
	OPENBSD_5_6:1.2.0.2
	OPENBSD_5_6_BASE:1.2
	OPENBSD_5_5:1.1.0.2
	OPENBSD_5_5_BASE:1.1;
locks; strict;
comment	@ * @;


1.4
date	2015.09.10.13.00.10;	author sthen;	state dead;
branches;
next	1.3;
commitid	vagA8GU9rOjsV69E;

1.3
date	2014.11.17.14.02.19;	author ajacoutot;	state Exp;
branches;
next	1.2;
commitid	lrMc8yxomUudhiIj;

1.2
date	2014.07.09.14.19.58;	author ajacoutot;	state Exp;
branches;
next	1.1;
commitid	bdiG2oZmTsT2fhbi;

1.1
date	2013.11.15.19.09.17;	author sthen;	state Exp;
branches;
next	;


desc
@@


1.4
log
@Currently Squid 3.5 is considered stable. Rather than moving snapshot/*
to stable/*, just drop back to a single version of squid in ports.
@
text
@#!/bin/sh
#
# $OpenBSD: squid.rc,v 1.3 2014/11/17 14:02:19 ajacoutot Exp $

daemon="${TRUEPREFIX}/sbin/squid"
daemon_timeout="35"

. /etc/rc.d/rc.subr

rc_stop() {
	${daemon} ${daemon_flags} -k shutdown
}

rc_reload() {
	${daemon} ${daemon_flags} -k reconfigure
}

rc_cmd $1
@


1.3
log
@Bump daemon_timeout from 31 to 35. Default squid shutdown timeout is 30s,
but we may be a bit slow, so give us some room.

ok sthen@@ (maintainer)
@
text
@d3 1
a3 1
# $OpenBSD: squid.rc,v 1.2 2014/07/09 14:19:58 ajacoutot Exp $
@


1.2
log
@Drop local shutdown patch and set daemon_timeout="31".

prodded by sthen@@
@
text
@d3 1
a3 1
# $OpenBSD: squid.rc,v 1.1 2013/11/15 19:09:17 sthen Exp $
d6 1
a6 1
daemon_timeout="31"
@


1.1
log
@Update to squid-3.3.10 and reorganise the squid ports; retire squid-2.7
(it was kept separately because 3.x didn't build on some arch due to an
incomplete atomic ops autoconf check which has since been repaired) and
replace with 3.4 (devel version). Enable additional helpers and turn on
support for rock cache dirs.
@
text
@d3 1
a3 1
# $OpenBSD: squid.rc,v 1.5 2011/05/18 14:47:55 sthen Exp $
d6 1
@

