# cvs/ports/devel/py-buildbot/pkg/buildbot.rc,v
head	1.4;
access;
symbols
	OPENBSD_6_1:1.4.0.4
	OPENBSD_6_1_BASE:1.4
	OPENBSD_6_0:1.4.0.2
	OPENBSD_6_0_BASE:1.4
	OPENBSD_5_9:1.3.0.2
	OPENBSD_5_9_BASE:1.3
	OPENBSD_5_8:1.3.0.4
	OPENBSD_5_8_BASE:1.3
	OPENBSD_5_7:1.2.0.2
	OPENBSD_5_7_BASE:1.2
	OPENBSD_5_6:1.1.0.14
	OPENBSD_5_6_BASE:1.1
	OPENBSD_5_5:1.1.0.12
	OPENBSD_5_5_BASE:1.1
	OPENBSD_5_4:1.1.0.10
	OPENBSD_5_4_BASE:1.1
	OPENBSD_5_3:1.1.0.8
	OPENBSD_5_3_BASE:1.1
	OPENBSD_5_2:1.1.0.6
	OPENBSD_5_2_BASE:1.1
	OPENBSD_5_1_BASE:1.1
	OPENBSD_5_1:1.1.0.4
	OPENBSD_5_0:1.1.0.2
	OPENBSD_5_0_BASE:1.1;
locks; strict;
comment	@ * @;


1.4
date	2016.05.15.09.04.09;	author landry;	state Exp;
branches;
next	1.3;
commitid	HmhFV2tLIrcwbOh5;

1.3
date	2015.07.16.00.09.32;	author ajacoutot;	state Exp;
branches;
next	1.2;
commitid	cZwJr0e7hjt7JI6F;

1.2
date	2015.01.05.20.12.30;	author landry;	state Exp;
branches;
next	1.1;
commitid	rCvN022XSQELKOpg;

1.1
date	2011.04.21.10.20.43;	author eric;	state Exp;
branches;
next	;


desc
@@


1.4
log
@Fix pexp in rc script to correctly detect status
ok/with ajacoutot@@
@
text
@#!/bin/sh
#
# $OpenBSD: buildbot.rc,v 1.3 2015/07/16 00:09:32 ajacoutot Exp $

daemon="${TRUEPREFIX}/bin/buildbot"
daemon_flags="${HOMEDIR}"
daemon_user="_buildbot"

. /etc/rc.d/rc.subr

pexp="${MODPY_BIN} .*logfile=twistd.log --python=buildbot.tac"

rc_start() {
	${rcexec} "${daemon} start ${daemon_flags}"
}

rc_reload() {
	${daemon} reconfig ${daemon_flags}
}

rc_stop() {
	${daemon} stop ${daemon_flags}
}

rc_cmd $1
@


1.3
log
@Hopefull fix pexp; this port is curently broken at runtime.
@
text
@d3 1
a3 1
# $OpenBSD: buildbot.rc,v 1.2 2015/01/05 20:12:30 landry Exp $
d11 1
a11 1
pexp="${MODPY_BIN} logfile=twistd.log --python=buildbot.tac"
@


1.2
log
@Fix pexp .. for some reason, the full commandline is not the same
anymore, and this breaks proper rc.d usage. Not ideal since this
probably prevents running several distinct instances, will be revisited
when i look at the 0.8.10 update.
@
text
@d3 1
a3 1
# $OpenBSD: buildbot.rc,v 1.1 2011/04/21 10:20:43 eric Exp $
d11 1
a11 1
pexp="${MODPY_BIN}.*logfile=twistd.log --python=buildbot.tac"
@


1.1
log
@- update to 0.8.3p1
- add a new _buildbot user
- add a rc script
- add a README
- update HOMEPAGE
- change MAINTAINER (per his request)

from Piotr Sikora (new MAINTAINER)

ok ajacoutot@@
@
text
@d3 1
a3 1
# $OpenBSD$
d11 1
a11 1
pexp="${MODPY_BIN} ${daemon} start ${daemon_flags}"
@

