# cvs/ports/www/transproxy/pkg/Attic/tproxy.rc,v
head	1.2;
access;
symbols
	OPENBSD_5_0:1.1.0.2
	OPENBSD_5_0_BASE:1.1;
locks; strict;
comment	@ * @;


1.2
date	2011.09.10.06.16.51;	author ajacoutot;	state dead;
branches;
next	1.1;

1.1
date	2011.03.07.07.25.40;	author ajacoutot;	state Exp;
branches;
next	;


desc
@@


1.2
log
@Remove tansproxy. It's not needed anymore (PF "route-to" works perfect
for that purpose).

ok sthen@@
@
text
@#!/bin/sh
#
# $OpenBSD: tproxy.rc,v 1.1 2011/03/07 07:25:40 ajacoutot Exp $

daemon="${TRUEPREFIX}/sbin/tproxy"
daemon_flags="-p -s 3128 -r _tproxy 127.0.0.1 3128"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
@


1.1
log
@Long overdue update to transproxy-1.6 (fixes a problem with posting
large amount of data amongst other things).

Add an rc.d(8) script.
Add a _tproxy user/group and make tproxy runs as this user by default.
Remove jakob@@ from maintainer per his request.

While here, do some cleanings:
* clean-up indent
* don't overwrite CC, CFLAGS and PREFIX
* fix patch naming
* fix a wording in DESCR
* comment tproxyrun and tproxywatch, they are useless under OpenBSD
@
text
@d3 1
a3 1
# $OpenBSD$
@

