# cvs/ports/net/unbound/pkg/Attic/unbound.rc,v
head	1.4;
access;
symbols
	OPENBSD_5_5:1.3.0.6
	OPENBSD_5_5_BASE:1.3
	OPENBSD_5_4:1.3.0.4
	OPENBSD_5_4_BASE:1.3
	OPENBSD_5_3:1.3.0.2
	OPENBSD_5_3_BASE:1.3
	OPENBSD_5_2:1.2.0.4
	OPENBSD_5_2_BASE:1.2
	OPENBSD_5_1_BASE:1.2
	OPENBSD_5_1:1.2.0.2
	OPENBSD_5_0:1.1.0.2
	OPENBSD_5_0_BASE:1.1;
locks; strict;
comment	@ * @;


1.4
date	2014.07.22.09.52.11;	author ajacoutot;	state dead;
branches;
next	1.3;
commitid	G5ZVuvgFae854xjI;

1.3
date	2012.08.04.20.43.54;	author sthen;	state Exp;
branches;
next	1.2;

1.2
date	2012.01.09.11.55.10;	author sthen;	state Exp;
branches;
next	1.1;

1.1
date	2011.03.06.09.24.41;	author jakob;	state Exp;
branches;
next	;


desc
@@


1.4
log
@Remove net/unbound. It's been unlinked for a while and is now in base.

prodded by and ok sthen@@, ok jasper@@
@
text
@#!/bin/sh
#
# $OpenBSD: unbound.rc,v 1.3 2012/08/04 20:43:54 sthen Exp $

daemon="${TRUEPREFIX}/sbin/unbound-control"
daemon_flags="-c /var/unbound/etc/unbound.conf"

. /etc/rc.d/rc.subr

pexp="unbound${daemon_flags:+ ${daemon_flags}}"
rc_reload=NO

rc_pre() {
	if ! [[ -f /var/unbound/etc/unbound_server.pem ||
	    -f /var/unbound/etc/unbound_control.key ||
	    -f /var/unbound/etc/unbound_control.pem ]]; then
		unbound-control-setup >/dev/null 2>&1
	fi
}

rc_start() {
	${rcexec} "${daemon} ${daemon_flags} start"
}

rc_check() {
	${daemon} ${daemon_flags} status
}

rc_stop() {
	${daemon} ${daemon_flags} stop
}

rc_cmd $1
@


1.3
log
@update to unbound 1.4.18; includes a crash fix (assertion failure).

- rc.d script now generates the unbound-control keys if they don't exist
and the sample config file is patched to enable this, various rc.d/unbound
actions depend on this, pointed out/ok aja@@
@
text
@d3 1
a3 1
# $OpenBSD: unbound.rc,v 1.2 2012/01/09 11:55:10 sthen Exp $
@


1.2
log
@Adjust unbound.rc to use unbound-control to stop/start/check Unbound.
Previously if unbound-control was used to restart the daemon, further
use of /etc/rc.d/unbound would no longer match the process title.
Reported by camield@@, ok aja@@ jakob@@
@
text
@d3 1
a3 1
# $OpenBSD: unbound.rc,v 1.1 2011/03/06 09:24:41 jakob Exp $
d12 8
@


1.1
log
@upgrade to Unbound 1.4.8; mostly from sthen@@
@
text
@d3 1
a3 1
# $OpenBSD$
d5 2
a6 1
daemon="${TRUEPREFIX}/sbin/unbound"
d9 15
@

