# cvs/ports/net/cnupm/pkg/cnupm.rc,v
head	1.1;
access;
symbols
	OPENBSD_6_2:1.1.0.8
	OPENBSD_6_2_BASE:1.1
	OPENBSD_6_1:1.1.0.6
	OPENBSD_6_1_BASE:1.1
	OPENBSD_6_0:1.1.0.4
	OPENBSD_6_0_BASE:1.1
	OPENBSD_5_9:1.1.0.2
	OPENBSD_5_9_BASE:1.1;
locks; strict;
comment	@ * @;


1.1
date	2015.08.14.13.06.52;	author zhuk;	state Exp;
branches;
next	;
commitid	1Tqx6DtBB834JPPi;


desc
@@


1.1
log
@Add rcscript for cnupm.

Input and okay aja@@.
@
text
@#!/bin/sh
#
# $OpenBSD$

daemon="${TRUEPREFIX}/sbin/cnupm"

. /etc/rc.d/rc.subr

cnupm_interface="${daemon_flags##*(* )-i*( )}"
cnupm_interface="${cnupm_interface%% *}"
[ -n "${cnupm_interface}" ] || cnupm_interface=".*"
pexp="cnupm: collecting traffic on ${cnupm_interface}"

rc_cmd $1
@
