# cvs/ports/audio/rplay/pkg/Attic/PLIST,v
head	1.18;
access;
symbols
	OPENBSD_5_9:1.17.0.2
	OPENBSD_5_9_BASE:1.17
	OPENBSD_5_8:1.17.0.4
	OPENBSD_5_8_BASE:1.17
	OPENBSD_5_7:1.16.0.6
	OPENBSD_5_7_BASE:1.16
	OPENBSD_5_6:1.16.0.4
	OPENBSD_5_6_BASE:1.16
	OPENBSD_5_5:1.16.0.2
	OPENBSD_5_5_BASE:1.16
	OPENBSD_5_4:1.15.0.32
	OPENBSD_5_4_BASE:1.15
	OPENBSD_5_3:1.15.0.30
	OPENBSD_5_3_BASE:1.15
	OPENBSD_5_2:1.15.0.28
	OPENBSD_5_2_BASE:1.15
	OPENBSD_5_1_BASE:1.15
	OPENBSD_5_1:1.15.0.26
	OPENBSD_5_0:1.15.0.24
	OPENBSD_5_0_BASE:1.15
	OPENBSD_4_9:1.15.0.22
	OPENBSD_4_9_BASE:1.15
	OPENBSD_4_8:1.15.0.20
	OPENBSD_4_8_BASE:1.15
	OPENBSD_4_7:1.15.0.18
	OPENBSD_4_7_BASE:1.15
	OPENBSD_4_6:1.15.0.16
	OPENBSD_4_6_BASE:1.15
	OPENBSD_4_5:1.15.0.14
	OPENBSD_4_5_BASE:1.15
	OPENBSD_4_4:1.15.0.12
	OPENBSD_4_4_BASE:1.15
	OPENBSD_4_3:1.15.0.10
	OPENBSD_4_3_BASE:1.15
	OPENBSD_4_2:1.15.0.8
	OPENBSD_4_2_BASE:1.15
	OPENBSD_4_1:1.15.0.6
	OPENBSD_4_1_BASE:1.15
	OPENBSD_4_0:1.15.0.4
	OPENBSD_4_0_BASE:1.15
	OPENBSD_3_9:1.15.0.2
	OPENBSD_3_9_BASE:1.15
	OPENBSD_3_8:1.14.0.2
	OPENBSD_3_8_BASE:1.14
	OPENBSD_3_7:1.13.0.2
	OPENBSD_3_7_BASE:1.13
	OPENBSD_3_6:1.12.0.2
	OPENBSD_3_6_BASE:1.12
	OPENBSD_3_5:1.11.0.2
	OPENBSD_3_5_BASE:1.11
	OPENBSD_3_4:1.10.0.12
	OPENBSD_3_4_BASE:1.10
	OPENBSD_3_3:1.10.0.10
	OPENBSD_3_3_BASE:1.10
	OPENBSD_3_2:1.10.0.8
	OPENBSD_3_2_BASE:1.10
	OPENBSD_3_1:1.10.0.6
	OPENBSD_3_1_BASE:1.10
	OPENBSD_3_0:1.10.0.4
	OPENBSD_3_0_BASE:1.10
	OPENBSD_2_9_TRACKING_SWITCH:1.10
	OPENBSD_2_9:1.10.0.2
	OPENBSD_2_9_BASE:1.10
	OPENBSD_2_8:1.9.0.2
	OPENBSD_2_8_BASE:1.9
	OPENBSD_2_7:1.6.0.2
	OPENBSD_2_7_BASE:1.6
	OPENBSD_2_6:1.5.0.4
	OPENBSD_2_6_BASE:1.5
	OPENBSD_2_5:1.5.0.2
	OPENBSD_2_5_BASE:1.5
	OPENBSD_2_4:1.4.0.2
	OPENBSD_2_4_BASE:1.4
	OPENBSD_2_3:1.3.0.2
	OPENBSD_2_3_BASE:1.3
	rplay-3_2_0b6:1.1.1.1
	angelos:1.1.1;
locks; strict;
comment	@# @;


1.18
date	2016.03.20.20.42.47;	author naddy;	state dead;
branches;
next	1.17;
commitid	IQk7bVJh3Ohcqszp;

1.17
date	2015.03.16.18.07.37;	author naddy;	state Exp;
branches;
next	1.16;
commitid	IJdm3Zowcdkcxp5X;

1.16
date	2014.01.11.11.22.43;	author espie;	state Exp;
branches;
next	1.15;

1.15
date	2005.12.23.15.07.06;	author espie;	state Exp;
branches;
next	1.14;

1.14
date	2005.04.17.21.18.29;	author espie;	state Exp;
branches;
next	1.13;

1.13
date	2004.09.22.10.02.18;	author espie;	state Exp;
branches;
next	1.12;

1.12
date	2004.08.06.14.37.28;	author espie;	state Exp;
branches;
next	1.11;

1.11
date	2003.09.30.10.36.36;	author naddy;	state Exp;
branches;
next	1.10;

1.10
date	2001.02.06.01.32.12;	author kevlo;	state Exp;
branches;
next	1.9;

1.9
date	2000.09.21.08.41.01;	author brad;	state Exp;
branches;
next	1.8;

1.8
date	2000.09.21.07.47.07;	author brad;	state Exp;
branches;
next	1.7;

1.7
date	2000.06.10.17.39.21;	author espie;	state Exp;
branches;
next	1.6;

1.6
date	2000.04.08.16.39.14;	author espie;	state Exp;
branches
	1.6.2.1;
next	1.5;

1.5
date	98.12.22.05.24.13;	author form;	state Exp;
branches;
next	1.4;

1.4
date	98.06.26.03.55.57;	author angelos;	state Exp;
branches;
next	1.3;

1.3
date	98.04.10.01.32.06;	author angelos;	state Exp;
branches;
next	1.2;

1.2
date	97.12.05.02.22.16;	author angelos;	state Exp;
branches;
next	1.1;

1.1
date	97.12.05.02.18.06;	author angelos;	state Exp;
branches
	1.1.1.1;
next	;

1.1.1.1
date	97.12.05.02.18.06;	author angelos;	state Exp;
branches;
next	;

1.6.2.1
date	2000.09.15.04.42.55;	author marc;	state Exp;
branches;
next	;


desc
@@


1.18
log
@Remove obsolete, broken software from the 1990s:
fvwm95, a window manager with a Windows 95 look; ok jca@@ jung@@ mmcc@@
rplay, a forgotten network audio system; ok sthen@@ jca@@
@
text
@@@comment $OpenBSD: PLIST,v 1.17 2015/03/16 18:07:37 naddy Exp $
@@bin bin/rplay
@@bin bin/rptp
include/rplay.h
@@comment include/rxposix.h
@@info info/RPLAY.info
@@info info/RPTP.info
@@info info/librplay.info
@@info info/rplayd.info
lib/librplay.a
@@lib lib/librplay.so.${LIBrplay_VERSION}
@@comment lib/librx.a
@@man man/man1/rplay.1
@@man man/man1/rptp.1
@@man man/man5/rplay.conf.5
@@man man/man5/rplay.helpers.5
@@man man/man5/rplay.hosts.5
@@man man/man5/rplay.servers.5
@@man man/man8/rplayd.8
@@bin sbin/rplayd
@@mode 755
@@owner root
@@group wheel
share/examples/rplay/
@@sample ${SYSCONFDIR}/rplay/
@@mode 644
share/examples/rplay/rplay.conf
@@sample ${SYSCONFDIR}/rplay/rplay.conf
share/examples/rplay/rplay.helpers
@@sample ${SYSCONFDIR}/rplay/rplay.helpers
share/examples/rplay/rplay.hosts
@@sample ${SYSCONFDIR}/rplay/rplay.hosts
share/examples/rplay/rplay.servers
@@sample ${SYSCONFDIR}/rplay/rplay.servers
@


1.17
log
@fold lib-only PFRAG.shared into PLIST, regen PLIST, and bump
@
text
@d1 1
a1 1
@@comment $OpenBSD: PLIST,v 1.16 2014/01/11 11:22:43 espie Exp $
@


1.16
log
@scrape @@endfake. I haven't been able to finish this in... ten years.
time to let it go.
@
text
@d1 3
a3 2
@@comment $OpenBSD: PLIST,v 1.15 2005/12/23 15:07:06 espie Exp $
%%SHARED%%
a5 5
lib/librplay.a
@@comment lib/librx.a
bin/rplay
bin/rptp
@@comment include/rxposix.h
d10 2
d20 1
a20 1
sbin/rplayd
@


1.15
log
@switch to SHARED_LIBS + some old conflict prevention
@
text
@d1 1
a1 1
@@comment $OpenBSD: PLIST,v 1.14 2005/04/17 21:18:29 espie Exp $
a6 1
@@endfake
@


1.14
log
@endfake markers/no need to install rplay,jpeg,png to build afterstep.
@
text
@d1 1
a1 1
@@comment $OpenBSD: PLIST,v 1.13 2004/09/22 10:02:18 espie Exp $
d10 1
d15 1
d24 1
a26 1
@@mode 755
@


1.13
log
@ditch INSTALL/DEINSTALL: sample functionality needs help from a
sample configuration, there.
@
text
@d1 7
a7 1
@@comment $OpenBSD: PLIST,v 1.12 2004/08/06 14:37:28 espie Exp $
a9 2
include/rplay.h
@@comment include/rxposix.h
a13 2
lib/librplay.a
@@comment lib/librx.a
a21 1
%%SHARED%%
@


1.12
log
@switch to new style PLISTS, using @@lib, ditching DYNLIBDIR.
Replace two INSTALL scripts with @@sample.
Reorder plists to match auto-generation, add comments for files that
don't get installed.
@
text
@d1 1
a1 1
@@comment $OpenBSD: PLIST,v 1.11 2003/09/30 10:36:36 naddy Exp $
d21 14
a34 6
@@cwd ${SYSCONFDIR}
@@extra rplay/rplay.conf
@@extra rplay/rplay.helpers
@@extra rplay/rplay.hosts
@@extra rplay/rplay.servers
@@extra rplay/
@


1.11
log
@@@extra handling of configuration files; ok kevlo@@
@
text
@d1 1
a1 1
@@comment $OpenBSD: PLIST,v 1.10 2001/02/06 01:32:12 kevlo Exp $
d5 5
a9 12
@@unexec install-info --delete --info-dir=%D/info %D/info/RPLAY.info
@@unexec install-info --delete --info-dir=%D/info %D/info/RPTP.info
@@unexec install-info --delete --info-dir=%D/info %D/info/librplay.info
@@unexec install-info --delete --info-dir=%D/info %D/info/rplayd.info
info/RPLAY.info
info/RPTP.info
info/librplay.info
info/rplayd.info
@@exec install-info --info-dir=%D/info %D/info/RPLAY.info
@@exec install-info --info-dir=%D/info %D/info/RPTP.info
@@exec install-info --info-dir=%D/info %D/info/librplay.info
@@exec install-info --info-dir=%D/info %D/info/rplayd.info
d11 9
a20 8
man/man1/rplay.1
man/man1/rptp.1
man/man5/rplay.conf.5
man/man5/rplay.helpers.5
man/man5/rplay.hosts.5
man/man5/rplay.servers.5
man/man8/rplayd.8
sbin/rplayd
@


1.10
log
@- Upgrade to version 3.2.2
- Take over maintainership

ok maintainer <angelos@@>
@
text
@d1 1
a1 1
@@comment $OpenBSD: PLIST,v 1.9 2000/09/21 08:41:01 brad Exp $
d27 6
@


1.9
log
@- rplay.info -> rplayd.info
- use better filenames for patches and use unified diffs
- remove post-build target
- in cdrom.c, stop trying to use /dev/acd
@
text
@d1 1
a1 1
@@comment $OpenBSD: PLIST,v 1.8 2000/09/21 07:47:07 brad Exp $
d18 1
@


1.8
log
@- remove license type and post-install target
- add HOMEPAGE to DESCR
- add proper INSTALL/DEINSTALL scripts
- silence configure script
@
text
@d1 1
a1 1
@@comment $OpenBSD: PLIST,v 1.7 2000/06/10 17:39:21 espie Exp $
d8 1
a8 1
@@unexec install-info --delete --info-dir=%D/info %D/info/rplay.info
d12 1
a12 1
info/rplay.info
d16 1
a16 1
@@exec install-info --info-dir=%D/info %D/info/rplay.info
@


1.7
log
@Various minor clean-ups:
- use new CONFIGURE_STYLE
- use PLIST/DESCR substitution to add HOMEPAGE and various fragments
- add systematic @@comment $OpenBSD$ to PLIST
- remove FAKE=Yes
@
text
@d1 4
a4 1
@@comment $OpenBSD$
a8 3
bin/rplay
bin/rptp
include/rplay.h
d13 4
a25 9
@@exec if [ ! -d /etc/rplay ]; then echo "Installing default configs..."; mkdir -p /etc/rplay; echo localhost >/etc/rplay/rplay.hosts; touch /etc/rplay/rplay.{conf,helpers,servers}; fi
@@unexec if [ -d /etc/rplay ]; then echo "***"; echo "*** Finished uninstalling rplay."; echo "*** Configuration files stored in /etc/rplay have not been removed."; echo "*** You may wish to remove them manually."; fi
@@exec install-info --info-dir=%D/info %D/info/RPLAY.info
@@exec install-info --info-dir=%D/info %D/info/RPTP.info
@@exec install-info --info-dir=%D/info %D/info/librplay.info
@@exec install-info --info-dir=%D/info %D/info/rplay.info
@@comment requires XForms
@@comment bin/xrplay
@@comment man/man1/xrplay.1
@


1.6
log
@Fake.
Install texinfo documentation properly.
@
text
@d1 1
@


1.6.2.1
log
@bring cut-off audio into 2.7 branch
@
text
@a0 1
@@comment $OpenBSD: PLIST,v 1.7 2000/06/10 17:39:21 espie Exp $
@


1.5
log
@/usr/local/etc -> /etc/rplay; OpenBSD porting policy.
btw, in 2.4 /usr/local/etc doesn't exists at all...
@
text
@d1 4
a6 3
@@comment requires XForms
@@comment bin/xrplay
sbin/rplayd
d8 2
a9 1
lib/librplay.a
d12 1
a12 2
info/RPLAY.info
info/RPTP.info
a13 1
@@comment man/man1/xrplay.1
d20 1
d23 7
@


1.4
log
@update to version 3.3.0
@
text
@a11 1
etc/rplay.hosts
d13 1
a13 2
man/man1/rplayd.8
man/man1/xrplay.1
d19 3
a21 1
etc/rc.d/rplayd.sh
@


1.3
log
@Updated, version 3.3.0a1.
@
text
@d12 1
d14 1
a14 1
man/man1/rplayd.1
d18 1
d21 1
@


1.2
log
@Oops.
@
text
@a0 2
etc/rc.d/rplayd.sh
etc/rplay.hosts
@


1.1
log
@Initial revision
@
text
@d14 7
a20 7
man/man1/rplay.1.gz
man/man1/rplayd.1.gz
man/man1/xrplay.1.gz
man/man1/rptp.1.gz
man/man5/rplay.conf.5.gz
man/man5/rplay.hosts.5.gz
man/man5/rplay.servers.5.gz
@


1.1.1.1
log
@Add rplay port.
@
text
@@
