# cvs/ports/sysutils/sysmon/pkg/sysmond.rc,v
head	1.2;
access;
symbols
	OPENBSD_6_2:1.2.0.4
	OPENBSD_6_2_BASE:1.2
	OPENBSD_6_1:1.2.0.2
	OPENBSD_6_1_BASE:1.2
	sthen_20160810:1.1.1.1
	sthen:1.1.1;
locks; strict;
comment	@ * @;


1.2
date	2016.08.20.16.41.47;	author ajacoutot;	state Exp;
branches;
next	1.1;
commitid	rJvc4h6l3ZXI8qxT;

1.1
date	2016.08.10.20.52.47;	author sthen;	state Exp;
branches
	1.1.1.1;
next	;
commitid	sRwqkOa1067y4dLR;

1.1.1.1
date	2016.08.10.20.52.47;	author sthen;	state Exp;
branches;
next	;
commitid	sRwqkOa1067y4dLR;


desc
@@


1.2
log
@Fix rc_reload() so that it works with daemon_rtable.

ok sthen@@
@
text
@#!/bin/sh
#
# $OpenBSD: sysmond.rc,v 1.1.1.1 2016/08/10 20:52:47 sthen Exp $

daemon="${TRUEPREFIX}/bin/sysmond"
daemon_flags="-f ${SYSCONFDIR}/sysmon.conf"
daemon_timeout="60"

. /etc/rc.d/rc.subr

rc_reload() {
	${daemon} ${daemon_flags} -t && \
		pkill -HUP -T "${daemon_rtable}" -xf "${pexp}"
}

rc_cmd $1
@


1.1
log
@Initial revision
@
text
@d3 1
a3 1
# $OpenBSD: sshd,v 1.3 2015/07/15 13:50:54 ajacoutot Exp $
d12 2
a13 1
	${daemon} ${daemon_flags} -t && pkill -HUP -xf "${pexp}"
@


1.1.1.1
log
@import sysutils/sysmon, from George Rosamond, ok landry

Sysmon is a network monitoring tool designed to provide high performance and
accurate network monitoring.

This tool is available in the public domain for anyone to use it that is
interested. It provides better performance and checking capabilities than other
tools such as Rover, Nocmon (not this: Nocmonitor), Whatsup, Big Brother, and
other such tools.

Configuration is simple and hierarchical, and can easily integrate email alerts.

Currently supported protocols include SMTP, IMAP, HTTP, TCP, UDP, NNTP, and
PING. Sysmon also provides a simple web output of monitored hosts.
@
text
@@
