# cvs/ports/www/nginx/stable/pkg/Attic/PLIST,v
head	1.2;
access;
symbols
	OPENBSD_5_6:1.1.0.8
	OPENBSD_5_6_BASE:1.1
	OPENBSD_5_5:1.1.0.6
	OPENBSD_5_5_BASE:1.1
	OPENBSD_5_4:1.1.0.4
	OPENBSD_5_4_BASE:1.1
	OPENBSD_5_3:1.1.0.2
	OPENBSD_5_3_BASE:1.1;
locks; strict;
comment	@# @;


1.2
date	2014.08.26.19.43.59;	author robert;	state dead;
branches;
next	1.1;
commitid	Y2vy3KpSALTztMGt;

1.1
date	2012.09.11.16.01.45;	author sthen;	state Exp;
branches;
next	;


desc
@@


1.2
log
@- remove the devel nginx packages
- apply chroot and other patches from the base system
- update to 1.7.4
@
text
@@@comment $OpenBSD: PLIST,v 1.1 2012/09/11 16:01:45 sthen Exp $
%%passenger%%
!%%passenger%%
@@man man/man8/nginx.8
@@bin sbin/nginx
@@mode 700
@@owner www
@@group www
@@sample ${NGINX_DIR}/tmp/
@@mode
@@owner
@@group
share/doc/pkg-readmes/${FULLPKGNAME}
share/nginx/
share/nginx/fastcgi_params
share/nginx/html/
@@sample ${NGINX_DIR}/html/
share/nginx/html/50x.html
@@sample ${NGINX_DIR}/html/50x.html
share/nginx/html/index.html
@@sample ${NGINX_DIR}/html/index.html
share/nginx/koi-utf
share/nginx/koi-win
share/nginx/mime.types
share/nginx/nginx.conf
share/nginx/scgi_params
share/nginx/uwsgi_params
share/nginx/win-utf
@@rcscript ${RCDIR}/enginx
@


1.1
log
@Update nginx stable in ports to 1.2.x and add devel (1.3.x), this is
separate from the version in the base OS which deliberately does not include
the mail proxy and additional modules. rc.d/enginx script is added (not
rc.d/nginx to avoid conflicting with base). Lua module is also supported.
Most of the work done by william@@, with a few tweaks by me and ok william.

README could use more work and as pointed out by ajacoutot we need to
review permissions of the tmp directory, but the diff is already long enough.
@@pkgpath markers are set to update from the 1.0.x packages to 1.2.x (with or
without the passenger flavour as appropriate).
@
text
@d1 1
a1 1
@@comment $OpenBSD: PLIST,v 1.9 2012/02/19 15:02:08 robert Exp $
@

