# cvs/ports/net/aircontrol/pkg/aircontrol.rc,v
head	1.7;
access;
symbols
	OPENBSD_6_2:1.7.0.10
	OPENBSD_6_2_BASE:1.7
	OPENBSD_6_1:1.7.0.8
	OPENBSD_6_1_BASE:1.7
	OPENBSD_6_0:1.7.0.6
	OPENBSD_6_0_BASE:1.7
	OPENBSD_5_9:1.7.0.2
	OPENBSD_5_9_BASE:1.7
	OPENBSD_5_8:1.7.0.4
	OPENBSD_5_8_BASE:1.7
	OPENBSD_5_7:1.6.0.16
	OPENBSD_5_7_BASE:1.6
	OPENBSD_5_6:1.6.0.14
	OPENBSD_5_6_BASE:1.6
	OPENBSD_5_5:1.6.0.12
	OPENBSD_5_5_BASE:1.6
	OPENBSD_5_4:1.6.0.10
	OPENBSD_5_4_BASE:1.6
	OPENBSD_5_3:1.6.0.8
	OPENBSD_5_3_BASE:1.6
	OPENBSD_5_2:1.6.0.6
	OPENBSD_5_2_BASE:1.6
	OPENBSD_5_1_BASE:1.6
	OPENBSD_5_1:1.6.0.4
	OPENBSD_5_0:1.6.0.2
	OPENBSD_5_0_BASE:1.6;
locks; strict;
comment	@ * @;


1.7
date	2015.07.15.19.49.41;	author ajacoutot;	state Exp;
branches;
next	1.6;
commitid	5Lv2La2t7Jf9rALu;

1.6
date	2011.07.01.10.54.53;	author sthen;	state Exp;
branches;
next	1.5;

1.5
date	2011.04.18.16.16.47;	author ajacoutot;	state Exp;
branches;
next	1.4;

1.4
date	2011.03.28.20.55.12;	author sthen;	state Exp;
branches;
next	1.3;

1.3
date	2011.03.17.16.49.27;	author ajacoutot;	state Exp;
branches;
next	1.2;

1.2
date	2011.03.09.18.51.48;	author ajacoutot;	state Exp;
branches;
next	1.1;

1.1
date	2011.03.03.13.20.43;	author sthen;	state Exp;
branches;
next	;


desc
@@


1.7
log
@Fix pexp.
@
text
@#!/bin/sh
#
# $OpenBSD: aircontrol.rc,v 1.6 2011/07/01 10:54:53 sthen Exp $

daemon="${TRUEPREFIX}/bin/catalina.sh"
daemon_flags="-Xmx768M -server"
daemon_user="_aircontrol"

. /etc/rc.d/rc.subr

pexp="$(${LOCALBASE}/bin/javaPathHelper -c aircontrol) .*-Dcatalina.base=/var/aircontrol"

rc_reload=NO

rc_start() {
        ${rcexec} \
	    "CATALINA_OPTS=\"${daemon_flags}\" \
	    JAVA_HOME=\"$(${LOCALBASE}/bin/javaPathHelper -h aircontrol)\" ${daemon} start"
}

rc_stop() {
	CATALINA_OPTS="${daemon_flags}" \
	    JAVA_HOME="$(${LOCALBASE}/bin/javaPathHelper -h aircontrol)" \
	    ${daemon} stop
}

rc_cmd $1
@


1.6
log
@fix rc_stop, pointed out by dcoppa@@
@
text
@d3 1
a3 1
# $OpenBSD: aircontrol.rc,v 1.5 2011/04/18 16:16:47 ajacoutot Exp $
d11 1
a11 1
pexp="$(${LOCALBASE}/bin/javaPathHelper -h aircontrol).*-Dcatalina.base=/var/aircontrol"
@


1.5
log
@Unless there's a very good reason for it, 'rcexec' should only be used in
rc_start, so remove from the local rc_stop.

ok sthen@@
@
text
@d3 1
a3 1
# $OpenBSD: aircontrol.rc,v 1.4 2011/03/28 20:55:12 sthen Exp $
d22 2
a23 2
	CATALINA_OPTS=\"${daemon_flags}\" \
	    JAVA_HOME=\"$(${LOCALBASE}/bin/javaPathHelper -h aircontrol)\" \
@


1.4
log
@set an explicit path to javaPathHelper, don't rely on it being in PATH
@
text
@d3 1
a3 1
# $OpenBSD: aircontrol.rc,v 1.3 2011/03/17 16:49:27 ajacoutot Exp $
d22 1
a22 2
        ${rcexec} \
	    "CATALINA_OPTS=\"${daemon_flags}\" \
d24 1
a24 1
	    ${daemon} stop"
@


1.3
log
@Remove all redirections -- rc.subr(8) now takes care of it.
@
text
@d3 1
a3 1
# $OpenBSD: aircontrol.rc,v 1.2 2011/03/09 18:51:48 ajacoutot Exp $
d11 1
a11 1
pexp="$(javaPathHelper -h aircontrol).*-Dcatalina.base=/var/aircontrol"
d18 1
a18 1
	    JAVA_HOME=\"$(javaPathHelper -h aircontrol)\" ${daemon} start"
d24 1
a24 1
	    JAVA_HOME=\"$(javaPathHelper -h aircontrol)\" \
@


1.2
log
@Remove now uneeded redirection.
@
text
@d3 1
a3 1
# $OpenBSD: aircontrol.rc,v 1.1 2011/03/03 13:20:43 sthen Exp $
d25 1
a25 1
	    ${daemon} stop" 2>/dev/null
@


1.1
log
@- update airControl to 1.3.6-beta

- add an rc script, which looks good to ajacoutot@@ (and might have
parts people can borrow if they want to write an rc script for other
Java daemons)
@
text
@d3 1
a3 1
# $OpenBSD$
d18 1
a18 2
	    JAVA_HOME=\"$(javaPathHelper -h aircontrol)\" \
	    ${daemon} start" >/dev/null
d25 1
a25 1
	    ${daemon} stop" >/dev/null 2>&1
@

