# cvs/ports/www/puppet-dashboard/pkg/puppet_dashboard_workers.rc,v
head	1.5;
access;
symbols
	OPENBSD_6_1:1.5.0.4
	OPENBSD_6_1_BASE:1.5
	OPENBSD_6_0:1.5.0.2
	OPENBSD_6_0_BASE:1.5
	OPENBSD_5_9:1.4.0.2
	OPENBSD_5_9_BASE:1.4
	OPENBSD_5_8:1.3.0.4
	OPENBSD_5_8_BASE:1.3
	OPENBSD_5_7:1.2.0.2
	OPENBSD_5_7_BASE:1.2
	OPENBSD_5_6:1.1.0.6
	OPENBSD_5_6_BASE:1.1
	OPENBSD_5_5:1.1.0.4
	OPENBSD_5_5_BASE:1.1
	OPENBSD_5_4:1.1.0.2
	OPENBSD_5_4_BASE:1.1;
locks; strict;
comment	@ * @;


1.5
date	2016.07.02.10.50.34;	author sthen;	state Exp;
branches;
next	1.4;
commitid	VGhlJVOmVfejuae3;

1.4
date	2015.09.03.20.33.42;	author jasper;	state Exp;
branches;
next	1.3;
commitid	eMFLqSjUr2KIQIft;

1.3
date	2015.07.15.14.59.15;	author ajacoutot;	state Exp;
branches;
next	1.2;
commitid	3M8rWDdk3egW952j;

1.2
date	2014.09.21.16.29.56;	author sebastia;	state Exp;
branches;
next	1.1;
commitid	PukWaorlQsdCTdFt;

1.1
date	2013.03.18.14.25.41;	author ajacoutot;	state Exp;
branches;
next	;


desc
@@


1.5
log
@remove LD_PRELOAD=libpthread.so from puppet-dashboard, and remove sebastia@@
from maintainer as he isn't using it any more.
@
text
@#!/bin/sh
#
# $OpenBSD: puppet_dashboard_workers.rc,v 1.4 2015/09/03 20:33:42 jasper Exp $

daemon="./script/delayed_job"
daemon_flags="-p dashboard -n 4 -m"
daemon_user="_puppet-dashboard"

. /etc/rc.d/rc.subr

pexp="ruby${MODRUBY_BINREV}:.*delayed_job.*"
rc_reload=NO

rc_check() {
	pgrep -xf "${pexp}"
}

rc_start() {
	${rcexec} "cd ${INSTDIR} && \
		env RAILS_ENV=production ${daemon} ${daemon_flags} start"
}

rc_stop() {
	${rcexec} "cd ${INSTDIR} && \
		env RAILS_ENV=production ${daemon} ${daemon_flags} stop"
}

rc_cmd $1
@


1.4
log
@fix pexp to unbreak rc_check()

ok aja@@ sebastia@@ (MAINTAINER)
@
text
@d3 1
a3 1
# $OpenBSD: puppet_dashboard_workers.rc,v 1.3 2015/07/15 14:59:15 ajacoutot Exp $
d20 1
a20 1
		env RAILS_ENV=production LD_PRELOAD=libpthread.so ${daemon} ${daemon_flags} start"
@


1.3
log
@Adapr rc.d script after recent rc.subr change.
@
text
@d3 1
a3 1
# $OpenBSD: puppet_dashboard_workers.rc,v 1.2 2014/09/21 16:29:56 sebastia Exp $
d11 1
a11 1
pexp="ruby${MODRUBY_BINREV}:.*delayed_job"
@


1.2
log
@Fix, or more an ugly hack, after mysql -> mariadb switch,
use LD_PRELOAD=libpthread.so., since investigation how to
best upgrade, might take a bit longer.
cluebat to the source of problem and fix from jeremy@@

OK jasper@@
@
text
@d3 1
a3 1
# $OpenBSD: puppet_dashboard_workers.rc,v 1.1 2013/03/18 14:25:41 ajacoutot Exp $
d15 1
a15 1
	pgrep -f "^${pexp}"
@


1.1
log
@Add a separate rc.d(8) script for the workers so that we can easily
change the flags in rc.conf.local(5).
@
text
@d3 1
a3 1
# $OpenBSD: puppet_dashboard.rc,v 1.1 2013/03/13 10:16:13 ajacoutot Exp $
d20 1
a20 1
		env RAILS_ENV=production ${daemon} ${daemon_flags} start"
@

