# cvs/ports/net/kea/pkg/kea.rc,v
head	1.2;
access;
symbols
	OPENBSD_6_2:1.2.0.8
	OPENBSD_6_2_BASE:1.2
	OPENBSD_6_1:1.2.0.6
	OPENBSD_6_1_BASE:1.2
	OPENBSD_6_0:1.2.0.4
	OPENBSD_6_0_BASE:1.2
	OPENBSD_5_9:1.2.0.2
	OPENBSD_5_9_BASE:1.2
	sthen_20151223:1.1.1.1
	sthen:1.1.1;
locks; strict;
comment	@ * @;


1.2
date	2015.12.23.15.58.55;	author ajacoutot;	state Exp;
branches;
next	1.1;
commitid	6dGqhRvakBrGbBM8;

1.1
date	2015.12.23.13.40.53;	author sthen;	state Exp;
branches
	1.1.1.1;
next	;
commitid	eTumTGLuGSSayHFP;

1.1.1.1
date	2015.12.23.13.40.53;	author sthen;	state Exp;
branches;
next	;
commitid	eTumTGLuGSSayHFP;


desc
@@


1.2
log
@Tweak pexp so that rc_check() does not need to be overriden.

ok sthen@@
@
text
@#!/bin/sh
#
# $OpenBSD: kea.rc,v 1.1.1.1 2015/12/23 13:40:53 sthen Exp $

daemon="${TRUEPREFIX}/sbin/keactrl"

. /etc/rc.d/rc.subr

pexp="${TRUEPREFIX}/sbin/kea-dhcp(4|6|-ddns).*"

rc_pre() {
	if ! [ -d ${LOCALSTATEDIR}/run/kea ]; then
		install -d -o root -g _kea -m 775 ${LOCALSTATEDIR}/run/kea
	fi
}

rc_start() {
	${rcexec} "${daemon} start"
}

rc_reload() {
	${daemon} reload
}

rc_stop() {
	${daemon} stop
}

rc_cmd $1
@


1.1
log
@Initial revision
@
text
@d3 1
a3 1
# $OpenBSD$
d9 1
a9 1
pexp="${TRUEPREFIX}/sbin/kea-dhcp(4|6|-ddns)"
a22 4
}

rc_check() {
	pgrep -q -f "${pexp}"
@


1.1.1.1
log
@import ports/net/kea, from Patrik Lundin, ok landry

Kea is a new open source DHCPv4/DHCPv6 server being developed by Internet
Systems Consortium. The objective of this project is to provide a very
high-performance, extensible DHCP server engine for use by enterprises and
service providers, either as is or with extensions and modifications.

Kea provides DHCPv4 and DHCPv6 servers, a dynamic DNS update module, a portable
DHCP library, libdhcp++, and a DHCP benchmarking tool, perfdhcp.
@
text
@@
