# cvs/ports/net/knot/pkg/knot.rc,v
head	1.5;
access;
symbols
	OPENBSD_6_2:1.5.0.2
	OPENBSD_6_2_BASE:1.5
	OPENBSD_6_1:1.4.0.8
	OPENBSD_6_1_BASE:1.4
	OPENBSD_6_0:1.4.0.6
	OPENBSD_6_0_BASE:1.4
	OPENBSD_5_9:1.4.0.2
	OPENBSD_5_9_BASE:1.4
	OPENBSD_5_8:1.4.0.4
	OPENBSD_5_8_BASE:1.4
	jca_20150412:1.1.1.1
	jca:1.1.1;
locks; strict;
comment	@ * @;


1.5
date	2017.06.01.13.42.01;	author jca;	state Exp;
branches;
next	1.4;
commitid	3UifAiZVbL5o2fn0;

1.4
date	2015.04.14.06.38.23;	author ajacoutot;	state Exp;
branches;
next	1.3;
commitid	sUhQA76oPnqNasfZ;

1.3
date	2015.04.12.17.09.20;	author jca;	state Exp;
branches;
next	1.2;
commitid	t444sANETwwRgUPs;

1.2
date	2015.04.12.16.26.10;	author jca;	state Exp;
branches;
next	1.1;
commitid	2pCioLa67I2Ll0Dz;

1.1
date	2015.04.12.16.07.13;	author jca;	state Exp;
branches
	1.1.1.1;
next	;
commitid	CpTWhahIOKHdZQm7;

1.1.1.1
date	2015.04.12.16.07.13;	author jca;	state Exp;
branches;
next	;
commitid	CpTWhahIOKHdZQm7;


desc
@@


1.5
log
@Major update to knot-2.4.2

Upstream stopped supporting LTS 1.6.* releases.  Let's move to the 2.x
branch which only supports LMDB (patched to use MDB_WRITEMAP) and
uses __thread (thus, needs clang or lang/gcc).

The config format changed a bit, /usr/local/libexec/knot1to2 might help.

Most work done by Alexis Vachette, ok fcambus@@
@
text
@#!/bin/sh
#
# $OpenBSD: knot.rc,v 1.4 2015/04/14 06:38:23 ajacoutot Exp $

daemon="${PREFIX}/sbin/knotd -d"
knotc="${PREFIX}/sbin/knotc"
knotcheckconf="${knotc} conf-check"

. /etc/rc.d/rc.subr

rc_pre() {
	install -d -o _knot -g wheel -m 750 /var/run/knot/
	${knotcheckconf}
}

rc_check() {
	${knotc} status
}

rc_reload() {
	${knotcheckconf} && ${knotc} reload
}

rc_stop() {
	${knotc} stop
}

rc_cmd $1
@


1.4
log
@Sanitize rc.d script:
- no need to override rc_start()
- only checkconf when it makes sense

ok jca@@ (maintainer)
@
text
@d3 1
a3 1
# $OpenBSD: knot.rc,v 1.3 2015/04/12 17:09:20 jca Exp $
d7 1
a7 1
knotcheckconf="${knotc} checkconf"
@


1.3
log
@Trim uneeded "|| return 1" clauses from knot.rc

Spotted by ajacoutot@@
@
text
@d3 1
a3 1
# $OpenBSD: knot.rc,v 1.2 2015/04/12 16:26:10 jca Exp $
d7 1
d11 3
a13 4
rc_start() {
	install -d -o _knot -g wheel -m 750 ${LOCALSTATEDIR}/run/knot/
	${knotc} checkconf || return 1
	${daemon} ${daemon_flags}
d21 1
a21 1
	${knotc} reload
@


1.2
log
@More fixes before import.

- --without-stuff -> --disable-stuff where needed
- replace /var and /etc with LOCALSTATEDIR and SYSCONFDIR
- add "OpenBSD" cvs tag to knot.rc

The two latter nits were pointed out by portcheck(1).
@
text
@d3 1
a3 1
# $OpenBSD$
d17 1
a17 1
	${knotc} status || return 1
d21 1
a21 1
	${knotc} reload || return 1
d25 1
a25 1
	${knotc} stop || return 1
@


1.1
log
@Initial revision
@
text
@d2 2
d11 1
a11 1
	install -d -o _knot -g wheel -m 750 /var/run/knot/
@


1.1.1.1
log
@Import knot-1.6.2

Tests by Pierre Emeriaud (co-maintainer) and florian@@, ok sthen@@

DESCR:
Knot DNS is a high-performance authoritative-only DNS server which
supports all key features of the domain name system including zone
transfers and DNSSEC.  Among its goals is support for on-line addition
and removal of zones.

@
text
@@
