# cvs/src/etc/rc.d/Attic/ipropd_slave,v
head	1.2;
access;
symbols
	OPENBSD_5_5:1.1.0.4
	OPENBSD_5_5_BASE:1.1;
locks; strict;
comment	@# @;


1.2
date	2014.04.23.09.56.56;	author ajacoutot;	state dead;
branches;
next	1.1;

1.1
date	2013.08.09.16.24.54;	author ajacoutot;	state Exp;
branches;
next	;


desc
@@


1.2
log
@Remove krb5 bits from rc(8).

ok reyk@@
@
text
@#!/bin/sh
#
# $OpenBSD: ipropd_slave,v 1.1 2013/08/09 16:24:54 ajacoutot Exp $

daemon="/usr/libexec/ipropd-slave --detach"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
@


1.1
log
@Add rc.d(8) scripts for ipropd-master and ipropd-slave.

ok deraadt@@
@
text
@d3 1
a3 1
# $OpenBSD$
@

