# cvs/src/etc/rc.d/Attic/ospfd6,v
head	1.2;
access;
symbols;
locks; strict;
comment	@# @;


1.2
date	2011.07.17.00.25.58;	author david;	state dead;
branches;
next	1.1;

1.1
date	2011.07.08.22.20.07;	author deraadt;	state Exp;
branches;
next	;


desc
@@


1.2
log
@use correct spelling of ospf6d script; ok deraadt@@
@
text
@#!/bin/sh
#
# $OpenBSD: ospfd6,v 1.1 2011/07/08 22:20:07 deraadt Exp $

daemon="/usr/sbin/ospfd6"

. /etc/rc.d/rc.subr

pexp="ospf6d: parent.*"

rc_cmd $1
@


1.1
log
@re-add ospfd/ospf6d support which got lost, and build scripts for them
@
text
@d3 1
a3 1
# $OpenBSD: bgpd,v 1.1 2011/07/06 18:55:36 robert Exp $
@

