# cvs/src/etc/rc.d/Attic/popa3d,v
head	1.4;
access;
symbols
	OPENBSD_5_4:1.3.0.8
	OPENBSD_5_4_BASE:1.3
	OPENBSD_5_3:1.3.0.6
	OPENBSD_5_3_BASE:1.3
	OPENBSD_5_2:1.3.0.4
	OPENBSD_5_2_BASE:1.3
	OPENBSD_5_1_BASE:1.3
	OPENBSD_5_1:1.3.0.2;
locks; strict;
comment	@# @;


1.4
date	2013.12.15.17.40.45;	author tedu;	state dead;
branches;
next	1.3;

1.3
date	2011.12.21.09.26.29;	author ajacoutot;	state Exp;
branches;
next	1.2;

1.2
date	2011.09.19.20.42.24;	author ajacoutot;	state Exp;
branches;
next	1.1;

1.1
date	2011.09.16.10.11.20;	author robert;	state Exp;
branches;
next	;


desc
@@


1.4
log
@remove popa3d etc tendrils
@
text
@#!/bin/sh
#
# $OpenBSD: popa3d,v 1.3 2011/12/21 09:26:29 ajacoutot Exp $

daemon="/usr/sbin/popa3d -D"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
@


1.3
log
@Move mandatory flags from daemon_flags to daemon.

ok sthen@@
@
text
@d3 1
a3 1
# $OpenBSD: popa3d,v 1.2 2011/09/19 20:42:24 ajacoutot Exp $
@


1.2
log
@Set defaults daemon_flags (same as advertised in rc.conf.local(5) so
that we can use these scripts right away with /etc/rc.d/XXX -f without
having to manually set _flags.

ok robert@@
@
text
@d3 1
a3 1
# $OpenBSD: popa3d,v 1.1 2011/09/16 10:11:20 robert Exp $
d5 1
a5 2
daemon="/usr/sbin/popa3d"
daemon_flags="-D"
@


1.1
log
@Add a script for popa3d to support running it outside of inetd.
@
text
@d3 1
a3 1
# $OpenBSD$
d6 1
@

