# cvs/ports/sysutils/backuppc/pkg/wwbackuppc.rc,v
head	1.2;
access;
symbols
	OPENBSD_6_2:1.2.0.10
	OPENBSD_6_2_BASE:1.2
	OPENBSD_6_1:1.2.0.8
	OPENBSD_6_1_BASE:1.2
	OPENBSD_6_0:1.2.0.6
	OPENBSD_6_0_BASE:1.2
	OPENBSD_5_9:1.2.0.2
	OPENBSD_5_9_BASE:1.2
	OPENBSD_5_8:1.2.0.4
	OPENBSD_5_8_BASE:1.2
	OPENBSD_5_7:1.1.1.1.0.14
	OPENBSD_5_7_BASE:1.1.1.1
	OPENBSD_5_6:1.1.1.1.0.12
	OPENBSD_5_6_BASE:1.1.1.1
	OPENBSD_5_5:1.1.1.1.0.10
	OPENBSD_5_5_BASE:1.1.1.1
	OPENBSD_5_4:1.1.1.1.0.8
	OPENBSD_5_4_BASE:1.1.1.1
	OPENBSD_5_3:1.1.1.1.0.6
	OPENBSD_5_3_BASE:1.1.1.1
	OPENBSD_5_2:1.1.1.1.0.4
	OPENBSD_5_2_BASE:1.1.1.1
	OPENBSD_5_1_BASE:1.1.1.1
	OPENBSD_5_1:1.1.1.1.0.2
	landry_20110920:1.1.1.1
	landry:1.1.1;
locks; strict;
comment	@ * @;


1.2
date	2015.07.15.22.40.09;	author ajacoutot;	state Exp;
branches;
next	1.1;
commitid	cteXsgpPHydsJ1aS;

1.1
date	2011.09.20.11.36.55;	author landry;	state Exp;
branches
	1.1.1.1;
next	;

1.1.1.1
date	2011.09.20.11.36.55;	author landry;	state Exp;
branches;
next	;


desc
@@


1.2
log
@Fix pexp.
@
text
@#!/bin/sh
#
# $OpenBSD: wwbackuppc.rc,v 1.1.1.1 2011/09/20 11:36:55 landry Exp $

daemon="${LOCALBASE}/bin/spawn-fcgi"
daemon_flags="-u_backuppc -g_backuppc -p1028 ${TRUEPREFIX}/bin/BackupPC_Admin.fcgi"

. /etc/rc.d/rc.subr

pexp="/usr/bin/perl ${TRUEPREFIX}/bin/BackupPC_Admin.fcgi"
rc_reload=NO

rc_cmd $1
@


1.1
log
@Initial revision
@
text
@d3 1
a3 1
# $OpenBSD$
d10 1
a10 1
pexp="/usr/bin/perl.*BackupPC_Admin.fcgi"
@


1.1.1.1
log
@Import BackupPC 3.2.1.

BackupPC is a high-performance, enterprise-grade system for backing
up Unix, WinXX, and MacOS PCs and laptops to a server's disk.
BackupPC is highly configurable and easy to install and maintain.
It is written in Perl and extracts backup data via rsync, or tar
over ssh/rsh/nfs.

With tweaks from and ok aja@@
@
text
@@
