# cvs/ports/net/apache-activemq/pkg/activemq.rc,v
head	1.3;
access;
symbols
	OPENBSD_6_1:1.3.0.16
	OPENBSD_6_1_BASE:1.3
	OPENBSD_6_0:1.3.0.14
	OPENBSD_6_0_BASE:1.3
	OPENBSD_5_9:1.3.0.10
	OPENBSD_5_9_BASE:1.3
	OPENBSD_5_8:1.3.0.12
	OPENBSD_5_8_BASE:1.3
	OPENBSD_5_7:1.3.0.8
	OPENBSD_5_7_BASE:1.3
	OPENBSD_5_6:1.3.0.6
	OPENBSD_5_6_BASE:1.3
	OPENBSD_5_5:1.3.0.4
	OPENBSD_5_5_BASE:1.3
	OPENBSD_5_4:1.3.0.2
	OPENBSD_5_4_BASE:1.3
	OPENBSD_5_3:1.1.1.1.0.2
	OPENBSD_5_3_BASE:1.1.1.1
	jasper_20130502:1.1.1.1
	jasper:1.1.1;
locks; strict;
comment	@ * @;


1.3
date	2013.03.13.16.30.11;	author ajacoutot;	state Exp;
branches;
next	1.2;

1.2
date	2013.03.04.15.48.41;	author ajacoutot;	state Exp;
branches;
next	1.1;

1.1
date	2013.02.05.15.08.45;	author jasper;	state Exp;
branches
	1.1.1.1;
next	;

1.1.1.1
date	2013.02.05.15.08.45;	author jasper;	state Exp;
branches;
next	;


desc
@@


1.3
log
@Revamp that port:
* don't use an /etc/default/activemq config file; it's too generic and
linuxish -- we'll use /etc/activemqrc
* no need to install the wrappers since we have rc scripts
* properly sample *all* configuration files
* run as unpriviledged user: _activemq by default
* remove useless README
* properly remove data files at @@extraunexec time
* proper rc_stop in the rc script
* take MAINTAINER

ok jasper@@ on a previous submission
@
text
@#!/bin/sh
#
# $OpenBSD: activemq.rc,v 1.2 2013/03/04 15:48:41 ajacoutot Exp $

daemon="${TRUEPREFIX}/activemq/bin/activemq"
daemon_user="_activemq"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_start() {
	${rcexec} "PATH=$PATH:$(dirname $(${LOCALBASE}/bin/javaPathHelper -c apache-activemq)) ${daemon} start"
}
 
rc_check() {
	PATH=$PATH:$(dirname $(${LOCALBASE}/bin/javaPathHelper -c apache-activemq)) ${daemon} status
}

rc_stop() {
	PATH=$PATH:$(dirname $(${LOCALBASE}/bin/javaPathHelper -c apache-activemq)) ${daemon} stop
}

rc_cmd $1
@


1.2
log
@Repair rc script which could never have worked...
@
text
@d3 1
a3 1
# $OpenBSD: activemq.rc,v 1.1.1.1 2013/02/05 15:08:45 jasper Exp $
d5 2
a6 1
daemon="${TRUEPREFIX}/bin/activemq"
d17 1
a17 1
	PATH=$PATH:$(dirname $(${LOCALBASE}/bin/javaPathHelper -c apache-activemq)) ${daemon} status   
d21 1
a21 2
	pid=$(rc_check | tail -1 | perl -pi -e 's,.*pid \x27(.*)\x27.*,$1,g')
	kill ${pid}
@


1.1
log
@Initial revision
@
text
@d3 1
a3 1
# $OpenBSD$
d12 1
a12 2
        ${rcexec} "env PATH=$PATH:$(dirname `javaPathHelper -c
apache-activemq`) ${daemon} start"
d16 1
a16 2
        env PATH=$PATH:$(dirname `javaPathHelper -c apache-activemq`)
${daemon} status   
d20 2
a21 2
        pid=$(rc_check | tail -1 | perl -pi -e 's,.*pid \x27(.*)\x27.*,$1,g')
        kill ${pid}
@


1.1.1.1
log
@import apache-activemq 5.7.0

Apache ActiveMQ is the most popular and powerful open source
messaging and Integration Patterns server.  Apache ActiveMQ is fast,
supports many Cross Language Clients and Protocols, comes with easy to
use Enterprise Integration Patterns and many advanced features while
fully supporting JMS 1.1 and J2EE 1.4.

rc.d help and ok aja@@
@
text
@@
