# cvs/ports/databases/redis/pkg/redis.rc,v
head	1.5;
access;
symbols
	OPENBSD_6_2:1.5.0.10
	OPENBSD_6_2_BASE:1.5
	OPENBSD_6_1:1.5.0.8
	OPENBSD_6_1_BASE:1.5
	OPENBSD_6_0:1.5.0.6
	OPENBSD_6_0_BASE:1.5
	OPENBSD_5_9:1.5.0.2
	OPENBSD_5_9_BASE:1.5
	OPENBSD_5_8:1.5.0.4
	OPENBSD_5_8_BASE:1.5
	OPENBSD_5_7:1.4.0.2
	OPENBSD_5_7_BASE:1.4
	OPENBSD_5_6:1.2.0.4
	OPENBSD_5_6_BASE:1.2
	OPENBSD_5_5:1.2.0.2
	OPENBSD_5_5_BASE:1.2
	OPENBSD_5_4:1.1.0.8
	OPENBSD_5_4_BASE:1.1
	OPENBSD_5_3:1.1.0.6
	OPENBSD_5_3_BASE:1.1
	OPENBSD_5_2:1.1.0.4
	OPENBSD_5_2_BASE:1.1
	OPENBSD_5_1_BASE:1.1
	OPENBSD_5_1:1.1.0.2;
locks; strict;
comment	@ * @;


1.5
date	2015.07.15.16.15.35;	author ajacoutot;	state Exp;
branches;
next	1.4;
commitid	xtIOzklA3EtpOLtM;

1.4
date	2014.11.11.10.12.54;	author dcoppa;	state Exp;
branches;
next	1.3;
commitid	10nBgTwawS8W03Ob;

1.3
date	2014.11.06.12.52.13;	author dcoppa;	state Exp;
branches;
next	1.2;
commitid	i62Xz51k9nFlCmfB;

1.2
date	2013.11.29.08.32.24;	author dcoppa;	state Exp;
branches;
next	1.1;

1.1
date	2011.11.29.16.07.23;	author dcoppa;	state Exp;
branches;
next	;


desc
@@


1.5
log
@Fix pexp after recent rc.subr change.
@
text
@#!/bin/sh
#
# $OpenBSD: redis.rc,v 1.4 2014/11/11 10:12:54 dcoppa Exp $

daemon="${TRUEPREFIX}/sbin/redis-server"
daemon_flags="${SYSCONFDIR}/redis/redis.conf"
daemon_user="_redis"

. /etc/rc.d/rc.subr

pexp="redis-server: ${daemon} .*"

rc_reload=NO

rc_pre() {
	/usr/bin/install -d -o _redis -g _redis -m 0750 /var/run/redis
}

rc_cmd $1
@


1.4
log
@
As pointed out by sthen@@, daemons should be in sbin not bin
@
text
@d3 1
a3 1
# $OpenBSD: redis.rc,v 1.3 2014/11/06 12:52:13 dcoppa Exp $
d11 1
a11 1
pexp="redis-server: ${daemon}.*"
@


1.3
log
@
Fix location/perms of the redis.sock file, for those who want to
use redis on a unix socket.

Amend redis.rc accordingly.

Also change the location of the pid file for consistency.

Discussed with, and ok landry@@
@
text
@d3 1
a3 1
# $OpenBSD: redis.rc,v 1.2 2013/11/29 08:32:24 dcoppa Exp $
d5 1
a5 1
daemon="${TRUEPREFIX}/bin/redis-server"
@


1.2
log
@Update to redis-2.8.1
@
text
@d3 1
a3 1
# $OpenBSD: redis.rc,v 1.1 2011/11/29 16:07:23 dcoppa Exp $
d14 4
@


1.1
log
@- Update to redis-2.4.4
- Enable backtracing support using libexecinfo's backtrace()
- Add a rc script

rc_script diff from william@@, who seems to have disappeared...
@
text
@d3 1
a3 1
# $OpenBSD$
d10 2
@

