# cvs/ports/sysutils/consul-template/pkg/consul_template.rc,v
head	1.1;
branch	1.1.1;
access;
symbols
	OPENBSD_6_1:1.1.1.1.0.2
	OPENBSD_6_1_BASE:1.1.1.1
	jasper_20160111:1.1.1.1
	jasper:1.1.1;
locks; strict;
comment	@ * @;


1.1
date	2016.11.01.16.51.30;	author jasper;	state Exp;
branches
	1.1.1.1;
next	;
commitid	t0gxZnOMxSUcRrnJ;

1.1.1.1
date	2016.11.01.16.51.30;	author jasper;	state Exp;
branches;
next	;
commitid	t0gxZnOMxSUcRrnJ;


desc
@@



1.1
log
@Initial revision
@
text
@#!/bin/sh
#
# $OpenBSD$

daemon="${PREFIX}/bin/consul-template"
daemon_flags="-config ${SYSCONFDIR}/consul-template.d/default.conf"
daemon_user="_consul-template"

. /etc/rc.d/rc.subr

rc_bg=YES

rc_stop() {
	pkill -INT -T "${daemon_rtable}" -xf "${pexp}"
}

rc_cmd $1
@


1.1.1.1
log
@import consul-template-0.16.0

This project provides a convenient way to populate values from Consul
into the file system using the consul-template daemon.

ok aja@@

@
text
@@
