# cvs/ports/sysutils/toad/pkg/Attic/toadd.rc,v
head	1.2;
access;
symbols
	OPENBSD_6_1:1.1.1.1.0.14
	OPENBSD_6_1_BASE:1.1.1.1
	OPENBSD_6_0:1.1.1.1.0.12
	OPENBSD_6_0_BASE:1.1.1.1
	OPENBSD_5_9:1.1.1.1.0.8
	OPENBSD_5_9_BASE:1.1.1.1
	OPENBSD_5_8:1.1.1.1.0.10
	OPENBSD_5_8_BASE:1.1.1.1
	OPENBSD_5_7:1.1.1.1.0.6
	OPENBSD_5_7_BASE:1.1.1.1
	OPENBSD_5_6:1.1.1.1.0.4
	OPENBSD_5_6_BASE:1.1.1.1
	OPENBSD_5_5:1.1.1.1.0.2
	OPENBSD_5_5_BASE:1.1.1.1
	ajacoutot_20131008:1.1.1.1
	ajacoutot:1.1.1;
locks; strict;
comment	@ * @;


1.2
date	2017.09.08.10.05.32;	author ajacoutot;	state dead;
branches;
next	1.1;
commitid	TyuAX2WwYnoakgeZ;

1.1
date	2013.10.08.17.27.01;	author ajacoutot;	state Exp;
branches
	1.1.1.1;
next	;

1.1.1.1
date	2013.10.08.17.27.01;	author ajacoutot;	state Exp;
branches;
next	;


desc
@@


1.2
log
@Remove toad, it requires kern.usermount -like functionality which is not
on the radar.
@
text
@#!/bin/sh
#
# $OpenBSD: toadd.rc,v 1.1 2013/10/08 17:27:01 ajacoutot Exp $

daemon="${TRUEPREFIX}/libexec/toadd"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
@


1.1
log
@Initial revision
@
text
@d3 1
a3 1
# $OpenBSD$
@


1.1.1.1
log
@Import toad-0.3.

toad(8) (Toad Opens All Devices) is a utility meant to be started from
the OpenBSD hotplugd(8) attach and detach scripts. It will try to mount
all partitions found on the device under /run/media/username/device.
Where username is the active user login name and device is the type of
the device: usb or cd, followed by its number (from 0 to 9).
This follows the udev hierarchy in Linux which allows interaction with
GLib/GIO's GUnixMount.

That means it interacts natively with file managers like Thunar and
Nautilus (where you can directly unmount/eject devices).

It comes with toadd(8), a small daemon that will detect the insertion of
a medium in the optical drives of the machine by periodically reading
their disklabel(8) (every 8 seconds) and call toad(8) with the attach or
detach argument respectively.

NOTE: this is somewhat still a WIP and a moving target for sure
ok jasper@@
@
text
@@
