# cvs/ports/www/mozilla-sync/full/pkg/mozilla_sync.rc,v
head	1.2;
access;
symbols
	OPENBSD_6_2:1.2.0.20
	OPENBSD_6_2_BASE:1.2
	OPENBSD_6_1:1.2.0.18
	OPENBSD_6_1_BASE:1.2
	OPENBSD_6_0:1.2.0.16
	OPENBSD_6_0_BASE:1.2
	OPENBSD_5_9:1.2.0.12
	OPENBSD_5_9_BASE:1.2
	OPENBSD_5_8:1.2.0.14
	OPENBSD_5_8_BASE:1.2
	OPENBSD_5_7:1.2.0.10
	OPENBSD_5_7_BASE:1.2
	OPENBSD_5_6:1.2.0.8
	OPENBSD_5_6_BASE:1.2
	OPENBSD_5_5:1.2.0.6
	OPENBSD_5_5_BASE:1.2
	OPENBSD_5_4:1.2.0.4
	OPENBSD_5_4_BASE:1.2
	OPENBSD_5_3:1.2.0.2
	OPENBSD_5_3_BASE:1.2
	landry_20120901:1.1.1.1
	landry:1.1.1;
locks; strict;
comment	@ * @;


1.2
date	2012.09.01.16.33.26;	author ajacoutot;	state Exp;
branches;
next	1.1;

1.1
date	2012.09.01.16.11.36;	author landry;	state Exp;
branches
	1.1.1.1;
next	;

1.1.1.1
date	2012.09.01.16.11.36;	author landry;	state Exp;
branches;
next	;


desc
@@


1.2
log
@rc.d script cleanup:
/etc -> SYSCONFDIR
also match on empty daemon_flags

ok landry@@
@
text
@#!/bin/sh
#
# $OpenBSD: mozilla_sync.rc,v 1.1.1.1 2012/09/01 16:11:36 landry Exp $

daemon="${LOCALBASE}/bin/paster"
daemon_flags="serve ${SYSCONFDIR}/mozilla-sync/production.ini --daemon --log-file /var/mozilla-sync/paster.log --pid-file /var/mozilla-sync/paster.pid"
daemon_user="_mozsync"

. /etc/rc.d/rc.subr

pexp="${MODPY_BIN} ${daemon}${daemon_flags:+ ${daemon_flags}}"
rc_reload=NO

rc_cmd $1
@


1.1
log
@Initial revision
@
text
@d3 1
a3 1
# $OpenBSD$
d6 1
a6 1
daemon_flags="serve /etc/mozilla-sync/production.ini --daemon --log-file /var/mozilla-sync/paster.log --pid-file /var/mozilla-sync/paster.pid"
d11 1
a11 2
pexp="${MODPY_BIN} ${daemon} ${daemon_flags}"

@


1.1.1.1
log
@Import Mozilla Sync Server.

Mozilla Sync server allows bookmarks, tabs, settings, addons, passwords
and history to be shared across several mozilla instances. Most work
done before previous lock, polished last month and finally here it is!
See https://support.mozilla.org/en-US/kb/how-do-i-set-up-firefox-sync
and the provided README for how to set it up properly.

ok sthen@@ on a previous version.
@
text
@@
