# cvs/ports/misc/memcached/pkg/memcached.rc,v
head	1.2;
access;
symbols
	OPENBSD_6_0:1.2.0.2
	OPENBSD_6_0_BASE:1.2
	OPENBSD_5_9:1.1.0.16
	OPENBSD_5_9_BASE:1.1
	OPENBSD_5_8:1.1.0.18
	OPENBSD_5_8_BASE:1.1
	OPENBSD_5_7:1.1.0.14
	OPENBSD_5_7_BASE:1.1
	OPENBSD_5_6:1.1.0.12
	OPENBSD_5_6_BASE:1.1
	OPENBSD_5_5:1.1.0.10
	OPENBSD_5_5_BASE:1.1
	OPENBSD_5_4:1.1.0.8
	OPENBSD_5_4_BASE:1.1
	OPENBSD_5_3:1.1.0.6
	OPENBSD_5_3_BASE:1.1
	OPENBSD_5_2:1.1.0.4
	OPENBSD_5_2_BASE:1.1
	OPENBSD_5_1_BASE:1.1
	OPENBSD_5_1:1.1.0.2;
locks; strict;
comment	@ * @;


1.2
date	2016.03.09.14.50.11;	author giovanni;	state Exp;
branches;
next	1.1;
commitid	kMrS9fQvG65cDUJa;

1.1
date	2011.12.25.14.37.28;	author sebastia;	state Exp;
branches;
next	;


desc
@@


1.2
log
@Add pledge(2) support
hints from semarie@@ and sthen@@
ok sthen@@
@
text
@#!/bin/sh
#
# $OpenBSD: memcached.rc,v 1.1 2011/12/25 14:37:28 sebastia Exp $

daemon="${TRUEPREFIX}/bin/memcached -d"
daemon_flags="-u _memcached"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
@


1.1
log
@update memcached 1.4.5 -> 1.4.10.
Added rc script.

tested with www/sogo on i386 by me, and with devel/subversion by stsp@@ on amd64

OK stsp@@
@
text
@d3 1
a3 1
# $OpenBSD$
d6 1
a6 1
daemon_flags="-u _memcached -P /var/run/memcached.pid"
@

