# cvs/ports/mail/zarafa/zarafa/pkg/Attic/zarafa-monitor.rc,v
head	1.2;
access;
symbols;
locks; strict;
comment	@ * @;


1.2
date	2010.10.28.20.46.02;	author ajacoutot;	state dead;
branches;
next	1.1;

1.1
date	2010.10.28.14.46.56;	author ajacoutot;	state Exp;
branches;
next	;


desc
@@


1.2
log
@After discussion with robert and ingo, we are renaming the rc scripts from
dash to underscore so that we are consitent with the foo_bar_flags
variable.
@
text
@#!/bin/sh

. /etc/rc.d/rc.subr

daemon="${TRUEPREFIX}/bin/zarafa-monitor"

rc_cmd $1
@


1.1
log
@Add zarafa daemons rc scripts.
Discussed with and ok robert@@
@
text
@@

