# cvs/ports/mail/zarafa/zarafa/pkg/Attic/zarafa_indexer.rc,v
head	1.2;
access;
symbols
	OPENBSD_5_2:1.1.0.6
	OPENBSD_5_2_BASE:1.1
	OPENBSD_5_1_BASE:1.1
	OPENBSD_5_1:1.1.0.4
	OPENBSD_5_0:1.1.0.2
	OPENBSD_5_0_BASE:1.1;
locks; strict;
comment	@ * @;


1.2
date	2012.08.30.14.28.22;	author ajacoutot;	state dead;
branches;
next	1.1;

1.1
date	2011.04.04.08.01.04;	author ajacoutot;	state Exp;
branches;
next	;


desc
@@


1.2
log
@Major update to zarafa-7.1.0.
Complete ChangeLog available at:
http://download.zarafa.com/community/final/7.1/final-changelog-7.1.txt
@
text
@#!/bin/sh
#
# $OpenBSD: zarafa_indexer.rc,v 1.1 2011/04/04 08:01:04 ajacoutot Exp $

daemon="${TRUEPREFIX}/bin/zarafa-indexer"

. /etc/rc.d/rc.subr

rc_cmd $1
@


1.1
log
@Add clucene support (i.e. enable zarafa-indexer). Note that the
zarafa-indexer is not used by default.
Tweak patches accordingly.
Use RLIMIT_DATA instead of RLIMIT_AS which we don't have.

Fix default client_update_path.
Remove a couple of patch entries that are not needed anymore.

Rewrite zarafa.rc and make it a wrapper script that calls the provided
rc.d(8) scripts. This is much more simple, makes much more sense and
allow for proper debugging as discussed with robert@@.
@
text
@d3 1
a3 1
# $OpenBSD$
@

