# cvs/ports/mail/zarafa/zarafa/pkg/Attic/zarafa_search.rc,v
head	1.3;
access;
symbols
	OPENBSD_5_8:1.2.0.4
	OPENBSD_5_8_BASE:1.2
	OPENBSD_5_7:1.1.0.10
	OPENBSD_5_7_BASE:1.1
	OPENBSD_5_6:1.1.0.8
	OPENBSD_5_6_BASE:1.1
	OPENBSD_5_5:1.1.0.6
	OPENBSD_5_5_BASE:1.1
	OPENBSD_5_4:1.1.0.4
	OPENBSD_5_4_BASE:1.1
	OPENBSD_5_3:1.1.0.2
	OPENBSD_5_3_BASE:1.1;
locks; strict;
comment	@ * @;


1.3
date	2015.09.28.19.05.24;	author robert;	state dead;
branches;
next	1.2;
commitid	0yNaOJnmdxPQsDK2;

1.2
date	2015.04.05.22.25.08;	author robert;	state Exp;
branches;
next	1.1;
commitid	j25w8UqbNvnTFexL;

1.1
date	2012.08.30.14.28.22;	author ajacoutot;	state Exp;
branches;
next	;


desc
@@


1.3
log
@remove zarafa_search because it works differently than before and an rc
script is not needed
@
text
@#!/bin/sh
#
# $OpenBSD: zarafa_search.rc,v 1.2 2015/04/05 22:25:08 robert Exp $

daemon="${TRUEPREFIX}/bin/zarafa-search"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_pre() {
	install -d -o _zabbix -g _zabbix /var/run/zabbix
}

rc_cmd $1
@


1.2
log
@update to 7.2.0
@
text
@d3 1
a3 1
# $OpenBSD: zarafa_search.rc,v 1.1 2012/08/30 14:28:22 ajacoutot Exp $
@


1.1
log
@Major update to zarafa-7.1.0.
Complete ChangeLog available at:
http://download.zarafa.com/community/final/7.1/final-changelog-7.1.txt
@
text
@d3 1
a3 1
# $OpenBSD: zarafa_indexer.rc,v 1.1 2011/04/04 08:01:04 ajacoutot Exp $
d10 4
@

