# cvs/ports/www/puppetboard/pkg/puppetboard.rc,v
head	1.6;
access;
symbols
	OPENBSD_6_0:1.5.0.6
	OPENBSD_6_0_BASE:1.5
	OPENBSD_5_9:1.5.0.2
	OPENBSD_5_9_BASE:1.5
	OPENBSD_5_8:1.5.0.4
	OPENBSD_5_8_BASE:1.5
	OPENBSD_5_7:1.3.0.2
	OPENBSD_5_7_BASE:1.3
	OPENBSD_5_6:1.2.0.2
	OPENBSD_5_6_BASE:1.2
	OPENBSD_5_5:1.1.0.2
	OPENBSD_5_5_BASE:1.1;
locks; strict;
comment	@ * @;


1.6
date	2016.11.02.09.01.13;	author danj;	state Exp;
branches;
next	1.5;
commitid	INr0zNQaJzHgaPPG;

1.5
date	2015.08.01.06.13.25;	author ajacoutot;	state Exp;
branches;
next	1.4;
commitid	0HCcD9d3DYSLZurb;

1.4
date	2015.07.15.22.24.51;	author ajacoutot;	state Exp;
branches;
next	1.3;
commitid	PXqj7aMR0Nz0ogPW;

1.3
date	2014.08.13.08.44.47;	author sebastia;	state Exp;
branches;
next	1.2;
commitid	Eq87BRqVCiEEqboH;

1.2
date	2014.06.17.11.55.11;	author jasper;	state Exp;
branches;
next	1.1;
commitid	ZktY818o0u15DWHF;

1.1
date	2013.12.27.11.06.33;	author jasper;	state Exp;
branches;
next	;


desc
@@


1.6
log
@Remove py3 flavor but make it default. For the rationale, see
https://marc.info/?l=openbsd-ports-cvs&m=147393010423200&w=2

While there, fix wrong substitution in the PLIST

test and OK sebastia@@ (who noticed the rc script needed to be updated)
@
text
@#!/bin/sh
#
# $OpenBSD: puppetboard.rc,v 1.5 2015/08/01 06:13:25 ajacoutot Exp $

daemon="puppetboard_wsgi"
daemon_flags="--workers 4"
daemon_user="_puppetboard"

. /etc/rc.d/rc.subr

pexp="${MODPY_BIN}.*${daemon_flags:+ ${daemon_flags}} ${daemon}:application"

rc_start() {
	${rcexec} "cd ${INSTDIR} && \
		${LOCALBASE}/bin/gunicorn${MODPY_BIN_SUFFIX} --daemon ${daemon_flags} ${daemon}:application"
}

rc_cmd $1
@


1.5
log
@Fix and extend pexp.

mismatch reported by sebastia@@
@
text
@d3 1
a3 1
# $OpenBSD: puppetboard.rc,v 1.4 2015/07/15 22:24:51 ajacoutot Exp $
d15 1
a15 1
		${LOCALBASE}/bin/gunicorn --daemon ${daemon_flags} ${daemon}:application"
@


1.4
log
@Fix pexp.
@
text
@d3 1
a3 1
# $OpenBSD: puppetboard.rc,v 1.3 2014/08/13 08:44:47 sebastia Exp $
d11 1
a11 1
pexp="${MODPY_BIN} .*${daemon}:application .*"
@


1.3
log
@puppetboard doesn't need to run as root, so add a user for it.

feedback and OK jasper@@ (MAINTAINER)
@
text
@d3 1
a3 1
# $OpenBSD: puppetboard.rc,v 1.2 2014/06/17 11:55:11 jasper Exp $
d11 1
a11 1
pexp="${MODPY_BIN}.*${daemon}:application"
@


1.2
log
@- fix ${pexp}

w/ aja@@
@
text
@d3 1
a3 1
# $OpenBSD: puppetboard.rc,v 1.1 2013/12/27 11:06:33 jasper Exp $
d7 1
@


1.1
log
@add rc script

ok aja@@
@
text
@d3 1
a3 1
# $OpenBSD$
d10 1
a10 1
pexp="python${MODPY_VERSION}:.*${daemon}:application"
@

