# cvs/ports/sysutils/ruby-puppet/pkg/Attic/puppetd.rc,v
head	1.8;
access;
symbols
	OPENBSD_5_7:1.7.0.4
	OPENBSD_5_7_BASE:1.7
	OPENBSD_5_6:1.7.0.2
	OPENBSD_5_6_BASE:1.7
	OPENBSD_5_5:1.6.0.2
	OPENBSD_5_5_BASE:1.6
	OPENBSD_5_4:1.5.0.2
	OPENBSD_5_4_BASE:1.5
	OPENBSD_5_3:1.3.0.2
	OPENBSD_5_3_BASE:1.3
	OPENBSD_5_2:1.2.0.4
	OPENBSD_5_2_BASE:1.2
	OPENBSD_5_1_BASE:1.2
	OPENBSD_5_1:1.2.0.2
	OPENBSD_5_0:1.1.0.2
	OPENBSD_5_0_BASE:1.1;
locks; strict;
comment	@ * @;


1.8
date	2015.04.04.09.28.45;	author jasper;	state dead;
branches;
next	1.7;
commitid	BK4zluMF2cw4Epmm;

1.7
date	2014.04.17.14.15.53;	author jasper;	state Exp;
branches;
next	1.6;

1.6
date	2013.11.12.08.28.29;	author robert;	state Exp;
branches;
next	1.5;

1.5
date	2013.03.13.14.35.39;	author ajacoutot;	state Exp;
branches;
next	1.4;

1.4
date	2013.03.08.16.12.48;	author ajacoutot;	state dead;
branches;
next	1.3;

1.3
date	2012.10.28.18.48.45;	author robert;	state Exp;
branches;
next	1.2;

1.2
date	2012.01.18.14.47.05;	author ajacoutot;	state Exp;
branches;
next	1.1;

1.1
date	2011.03.20.15.46.34;	author ajacoutot;	state Exp;
branches;
next	;


desc
@@


1.8
log
@rename rc scripts to drop the trailing 'd', so puppetd becomes just puppet.

this matches the upstream service names, as discussed with aja@@
@
text
@#!/bin/sh
#
# $OpenBSD: puppetd.rc,v 1.7 2014/04/17 14:15:53 jasper Exp $

daemon="${TRUEPREFIX}/bin/puppet agent"

. /etc/rc.d/rc.subr

pexp="${RUBY} ${daemon}${daemon_flags:+ ${daemon_flags}}$"
rc_reload=NO

rc_cmd $1
@


1.7
log
@- remove rc_stop() hack from puppetd, the issue has long since been resolved
- fixup license tag
@
text
@d3 1
a3 1
# $OpenBSD: puppetd.rc,v 1.6 2013/11/12 08:28:29 robert Exp $
@


1.6
log
@use a better regexp for the puppetd process to avoid killing every puppetd process; ok jasper@@
@
text
@d3 1
a3 1
# $OpenBSD: puppetd.rc,v 1.5 2013/03/13 14:35:39 ajacoutot Exp $
a10 6

# XXX can't stop within a fair amount of time when puppetmasterd is unreachable;
# related to: http://projects.puppetlabs.com/issues/11360 ?
rc_stop() {
	pkill -9 -f "^${pexp}"
}
@


1.5
log
@Factorize a bit more since both ports version share a lot.
While here, fix the rc script to start the master instead of the client
and vice-versa.

ok robert@@ (maintainer)
@
text
@d3 1
a3 1
# $OpenBSD: puppetd.rc,v 1.2 2013/03/13 09:56:19 ajacoutot Exp $
d9 1
a9 1
pexp="${RUBY} ${daemon}${daemon_flags:+ ${daemon_flags}}"
@


1.4
log
@Split ruby-puppet:
* move the current port to ruby-puppet/3
* add a new puppet-2.7.X port under ruby-puppet/2

This allows us to ship with 2 different versions of Puppet. Since the
Enterprise version is still running 2.7, we want to have a matching
client because running puppet3 against a puppet2 server is highly
discouraged and will mostly no work correctly.

ok jasper@@, robert@@ (maintainer)
@
text
@d3 1
a3 1
# $OpenBSD: puppetd.rc,v 1.3 2012/10/28 18:48:45 robert Exp $
d5 1
a5 1
daemon="/usr/local/bin/puppet agent"
d9 1
a9 1
pexp=".*ruby.* ${daemon}${daemon_flags:+ ${daemon_flags}}"
@


1.3
log
@update to 3.0.1
@
text
@d3 1
a3 1
# $OpenBSD: puppetd.rc,v 1.2 2012/01/18 14:47:05 ajacoutot Exp $
@


1.2
log
@XXX use pkill -9 in rc_stop for now; when puppetd is started then
stopped when puppetmasterd is not available it would fail to stop within
at least 60 seconds.

discussed with and ok robert@@ (maintainer)
@
text
@d3 1
a3 1
# $OpenBSD: puppetd.rc,v 1.1 2011/03/20 15:46:34 ajacoutot Exp $
d5 1
a5 1
daemon="${TRUEPREFIX}/sbin/puppetd"
@


1.1
log
@Add rc scripts for puppetd and puppetmasterd.
Set GPL version while here.

ok robert@@ (maintainer)
@
text
@d3 1
a3 1
# $OpenBSD$
d11 6
@

