# cvs/ports/net/openfire/pkg/openfire.rc,v
head	1.3;
access;
symbols
	OPENBSD_6_0:1.3.0.14
	OPENBSD_6_0_BASE:1.3
	OPENBSD_5_9:1.3.0.10
	OPENBSD_5_9_BASE:1.3
	OPENBSD_5_8:1.3.0.12
	OPENBSD_5_8_BASE:1.3
	OPENBSD_5_7:1.3.0.8
	OPENBSD_5_7_BASE:1.3
	OPENBSD_5_6:1.3.0.6
	OPENBSD_5_6_BASE:1.3
	OPENBSD_5_5:1.3.0.4
	OPENBSD_5_5_BASE:1.3
	OPENBSD_5_4:1.3.0.2
	OPENBSD_5_4_BASE:1.3
	OPENBSD_5_3:1.2.0.8
	OPENBSD_5_3_BASE:1.2
	OPENBSD_5_2:1.2.0.6
	OPENBSD_5_2_BASE:1.2
	OPENBSD_5_1_BASE:1.2
	OPENBSD_5_1:1.2.0.4
	OPENBSD_5_0:1.2.0.2
	OPENBSD_5_0_BASE:1.2;
locks; strict;
comment	@ * @;


1.3
date	2013.03.04.15.56.13;	author ajacoutot;	state Exp;
branches;
next	1.2;

1.2
date	2011.03.09.17.08.00;	author ajacoutot;	state Exp;
branches;
next	1.1;

1.1
date	2011.03.09.16.50.40;	author ian;	state Exp;
branches;
next	;


desc
@@


1.3
log
@Explicitely set path to javaPathHelper; we do not want to rely whether
LOCALBASE/bin is in the PATH or not and this is consistent with what
other java rc scripts do.

prodded by Marc Peters
@
text
@#!/bin/sh
#
# $OpenBSD: openfire.rc,v 1.2 2011/03/09 17:08:00 ajacoutot Exp $

OPENFIRE_HOME="${TRUEPREFIX}/openfire"
OPENFIRE_LIB="${OPENFIRE_HOME}/lib"

JAVACMD="$(${LOCALBASE}/bin/javaPathHelper -c openfire)"
JAVA_HOME="$(${LOCALBASE}/bin/javaPathHelper -h openfire)"

daemon="${JAVACMD}"
daemon_flags="-server -DopenfireHome=${OPENFIRE_HOME} -Dopenfire.lib.dir=${OPENFIRE_LIB} -jar ${OPENFIRE_LIB}/startup.jar"
daemon_user="_openfire"

. /etc/rc.d/rc.subr

rc_bg=YES
rc_reload=NO

rc_cmd $1
@


1.2
log
@I'm a consistency freak... no functionnal change it's only cosmetic.
@
text
@d3 1
a3 1
# $OpenBSD: openfire.rc,v 1.1 2011/03/09 16:50:40 ian Exp $
d8 2
a9 2
JAVACMD="$(javaPathHelper -c openfire)"
JAVA_HOME="$(javaPathHelper -h openfire)"
@


1.1
log
@Major update to 3.7.0. License changed to ASL2.0. Add rc.d script,
replacing openfirectl, using javaPathHelper. set MODJAVA_JRE =  Yes.
Upstream incorporated my previous patches, so rm them.  Remove MESSAGE
(missed during bulk conversion to pkg-readmes). Tweak README.

Requires a fairly current rc.subr - uses rc_bg !!!

Improved & OK sthen@@
@
text
@d2 2
a3 2

# $OpenBSD$
@

