# cvs/ports/sysutils/salt/pkg/salt_syndic.rc,v
head	1.2;
access;
symbols
	OPENBSD_6_1:1.2.0.12
	OPENBSD_6_1_BASE:1.2
	OPENBSD_6_0:1.2.0.10
	OPENBSD_6_0_BASE:1.2
	OPENBSD_5_9:1.2.0.6
	OPENBSD_5_9_BASE:1.2
	OPENBSD_5_8:1.2.0.8
	OPENBSD_5_8_BASE:1.2
	OPENBSD_5_7:1.2.0.4
	OPENBSD_5_7_BASE:1.2
	OPENBSD_5_6:1.2.0.2
	OPENBSD_5_6_BASE:1.2;
locks; strict;
comment	@ * @;


1.2
date	2014.06.17.12.30.39;	author ajacoutot;	state Exp;
branches;
next	1.1;
commitid	x81xXgb42w072Lr6;

1.1
date	2014.05.27.14.12.43;	author ajacoutot;	state Exp;
branches;
next	;


desc
@@


1.2
log
@The homedir of the user running salt must be outside of the /srv/salt
hierarchy because content of that path is available to all minions; so
move it to /var/salt .
Rename _salt-master to _salt for consistency.
discussed with "viq"

While here, drop rc_pre() from the minion and syndir rc.d(8) scripts,
they are not needed.
@
text
@#!/bin/sh
#
# $OpenBSD: salt_syndic.rc,v 1.1 2014/05/27 14:12:43 ajacoutot Exp $

daemon="${TRUEPREFIX}/bin/salt-syndic -d"
daemon_flags="-u _salt"

. /etc/rc.d/rc.subr

pexp="${MODPY_BIN} ${daemon}${daemon_flags:+ ${daemon_flags}}"
rc_reload=NO

rc_cmd $1
@


1.1
log
@Simplify this port not to deviate too much from upstream defaults.
Properly run as _salt-master without breaking.
Add an rc.d script for salt-syndic.
Drop useless README.
Remove all files when running pkg_delete -c
Install more config file examples.
Take maintainer.
@
text
@d3 1
a3 1
# $OpenBSD: salt_minion.rc,v 1.1.1.1 2013/08/10 13:39:27 landry Exp $
d6 1
a6 1
daemon_flags="-u _salt-master"
a8 4

rc_pre() {
	install -o _salt-master -d /var/run/salt/
}
@

