# cvs/ports/audio/mpd/pkg/mpd.rc,v
head	1.4;
access;
symbols
	OPENBSD_6_2:1.4.0.28
	OPENBSD_6_2_BASE:1.4
	OPENBSD_6_1:1.4.0.26
	OPENBSD_6_1_BASE:1.4
	OPENBSD_6_0:1.4.0.24
	OPENBSD_6_0_BASE:1.4
	OPENBSD_5_9:1.4.0.20
	OPENBSD_5_9_BASE:1.4
	OPENBSD_5_8:1.4.0.22
	OPENBSD_5_8_BASE:1.4
	OPENBSD_5_7:1.4.0.18
	OPENBSD_5_7_BASE:1.4
	OPENBSD_5_6:1.4.0.16
	OPENBSD_5_6_BASE:1.4
	OPENBSD_5_5:1.4.0.14
	OPENBSD_5_5_BASE:1.4
	OPENBSD_5_4:1.4.0.12
	OPENBSD_5_4_BASE:1.4
	OPENBSD_5_3:1.4.0.10
	OPENBSD_5_3_BASE:1.4
	OPENBSD_5_2:1.4.0.8
	OPENBSD_5_2_BASE:1.4
	OPENBSD_5_1_BASE:1.4
	OPENBSD_5_1:1.4.0.6
	OPENBSD_5_0:1.4.0.4
	OPENBSD_5_0_BASE:1.4
	OPENBSD_4_9:1.4.0.2
	OPENBSD_4_9_BASE:1.4;
locks; strict;
comment	@ * @;


1.4
date	2010.12.27.14.50.22;	author ajacoutot;	state Exp;
branches;
next	1.3;

1.3
date	2010.12.24.10.40.04;	author ajacoutot;	state Exp;
branches;
next	1.2;

1.2
date	2010.10.29.12.54.43;	author ajacoutot;	state Exp;
branches;
next	1.1;

1.1
date	2010.10.28.11.24.35;	author jasper;	state Exp;
branches;
next	;


desc
@@


1.4
log
@Simplify after recent rc.subr change.
The framework is now stable and we will start documenting it (at last).
@
text
@#!/bin/sh
#
# $OpenBSD: mpd.rc,v 1.3 2010/12/24 10:40:04 ajacoutot Exp $

daemon="${TRUEPREFIX}/sbin/mpd"

. /etc/rc.d/rc.subr

rc_cmd $1
@


1.3
log
@Cope with recent rc.subr changes.
@
text
@d3 1
a3 3
# $OpenBSD: mpd.rc,v 1.2 2010/10/29 12:54:43 ajacoutot Exp $

. /etc/rc.d/rc.subr
d7 1
a7 1
rc_conf
@


1.2
log
@Add RCS IDs to rc scripts.
@
text
@d3 1
a3 1
# $OpenBSD$
d8 2
@


1.1
log
@- add rc script for mpd

ok robert@@
@
text
@d2 2
@

