# cvs/ports/databases/mariadb/pkg/mysqld.rc,v
head	1.3;
access;
symbols
	OPENBSD_6_1:1.2.0.8
	OPENBSD_6_1_BASE:1.2
	OPENBSD_6_0:1.2.0.6
	OPENBSD_6_0_BASE:1.2
	OPENBSD_5_9:1.2.0.2
	OPENBSD_5_9_BASE:1.2
	OPENBSD_5_8:1.2.0.4
	OPENBSD_5_8_BASE:1.2
	OPENBSD_5_7:1.1.1.1.0.8
	OPENBSD_5_7_BASE:1.1.1.1
	OPENBSD_5_6:1.1.1.1.0.6
	OPENBSD_5_6_BASE:1.1.1.1
	OPENBSD_5_5:1.1.1.1.0.4
	OPENBSD_5_5_BASE:1.1.1.1
	OPENBSD_5_4:1.1.1.1.0.2
	OPENBSD_5_4_BASE:1.1.1.1
	brad_20130404:1.1.1.1
	brad:1.1.1;
locks; strict;
comment	@ * @;


1.3
date	2017.05.27.13.53.06;	author ajacoutot;	state Exp;
branches;
next	1.2;
commitid	SDpZdO8hfYIlMQ9j;

1.2
date	2015.07.15.15.52.29;	author ajacoutot;	state Exp;
branches;
next	1.1;
commitid	1heQdbgx3JmVR9Ch;

1.1
date	2013.04.04.18.58.52;	author brad;	state Exp;
branches
	1.1.1.1;
next	;

1.1.1.1
date	2013.04.04.18.58.52;	author brad;	state Exp;
branches;
next	;


desc
@@


1.3
log
@Cope with recent rc.subr(8) change (removal of rc_bg etc.)
@
text
@#!/bin/sh
#
# $OpenBSD: mysqld.rc,v 1.2 2015/07/15 15:52:29 ajacoutot Exp $

daemon="${TRUEPREFIX}/bin/mysqld_safe"

. /etc/rc.d/rc.subr

pexp="${TRUEPREFIX}/libexec/mysqld .*"
rc_reload=NO

rc_cmd $1
@


1.2
log
@Fix rc.d script after recent rc.subr change.
@
text
@d3 1
a3 1
# $OpenBSD: mysqld.rc,v 1.1.1.1 2013/04/04 18:58:52 brad Exp $
a9 1
rc_bg=YES
@


1.1
log
@Initial revision
@
text
@d3 1
a3 1
# $OpenBSD: mysqld.rc,v 1.12 2011/03/08 10:46:31 ajacoutot Exp $
d9 1
a9 1
pexp="${TRUEPREFIX}/libexec/mysqld"
@


1.1.1.1
log
@Import of MariaDB.

MariaDB is a multi-threaded SQL database with a command syntax very similar
to mSQL.

ok sthen@@ giovanni@@
@
text
@@
