# cvs/ports/net/gdnsd/pkg/gdnsd.rc,v
head	1.4;
access;
symbols
	OPENBSD_6_1:1.3.0.6
	OPENBSD_6_1_BASE:1.3
	OPENBSD_6_0:1.3.0.4
	OPENBSD_6_0_BASE:1.3
	OPENBSD_5_9:1.3.0.2
	OPENBSD_5_9_BASE:1.3
	OPENBSD_5_8:1.2.0.6
	OPENBSD_5_8_BASE:1.2
	OPENBSD_5_7:1.2.0.2
	OPENBSD_5_7_BASE:1.2
	jasper_20140310:1.1.1.1
	jasper:1.1.1;
locks; strict;
comment	@ * @;


1.4
date	2017.08.10.20.14.56;	author jung;	state Exp;
branches;
next	1.3;
commitid	eDkNoSM3Y9ec6X7r;

1.3
date	2016.02.06.07.45.08;	author ajacoutot;	state Exp;
branches;
next	1.2;
commitid	dWGqJ2BUwfG2Zfeq;

1.2
date	2014.10.03.11.53.23;	author jasper;	state Exp;
branches;
next	1.1;
commitid	YnCW657LSFnl923G;

1.1
date	2014.10.03.11.31.27;	author jasper;	state Exp;
branches
	1.1.1.1;
next	;
commitid	sOXaUGRoDecFdq3a;

1.1.1.1
date	2014.10.03.11.31.27;	author jasper;	state Exp;
branches;
next	;
commitid	sOXaUGRoDecFdq3a;


desc
@@


1.4
log
@remove reload and switch default state dir, also take maintainer

ok sthen
@
text
@#!/bin/sh
#
# $OpenBSD: gdnsd.rc,v 1.3 2016/02/06 07:45:08 ajacoutot Exp $

daemon="${PREFIX}/sbin/gdnsd"
daemon_flags="start"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_restart() {
	${daemon} restart
}

rc_stop() {
	${daemon} stop
}

rc_cmd $1
@


1.3
log
@Unbreak rc.d script.
reported by jung@@
@
text
@d3 1
a3 1
# $OpenBSD: gdnsd.rc,v 1.2 2014/10/03 11:53:23 jasper Exp $
d10 1
a10 3
rc_reload() {
	${daemon} reload
}
@


1.2
log
@rcexec should only really be used in rc_start, unless there's a good reason to
use it elsewhere

prompted by aja@@
ok aja@@ sthen@@
@
text
@d3 1
a3 1
# $OpenBSD: gdnsd.rc,v 1.1.1.1 2014/10/03 11:31:27 jasper Exp $
d6 1
a8 4

rc_start() {
	${rcexec} "${daemon} start"
}
@


1.1
log
@Initial revision
@
text
@d3 1
a3 1
# $OpenBSD$
d14 1
a14 1
	${rcexec} "${daemon} reload"
d18 1
a18 1
	${rcexec} "${daemon} restart"
d22 1
a22 1
	${rcexec} "${daemon} stop"
@


1.1.1.1
log
@import gdnsd-1.11.5

gdnsd is an Authoritative-only DNS server. The initial g stands for
Geographic, as gdnsd offers a plugin system for geographic (or other
sorts of) balancing, redirection, and service-state-conscious failover.
The plugin system can also do things like weighted address/cname
records.

feedback/ok sthen@@
@
text
@@
