# cvs/ports/x11/kde4/workspace/pkg/kdm.rc,v
head	1.9;
access;
symbols
	OPENBSD_6_1:1.8.0.6
	OPENBSD_6_1_BASE:1.8
	OPENBSD_6_0:1.8.0.4
	OPENBSD_6_0_BASE:1.8
	OPENBSD_5_9:1.8.0.2
	OPENBSD_5_9_BASE:1.8
	OPENBSD_5_8:1.4.0.4
	OPENBSD_5_8_BASE:1.4
	OPENBSD_5_7:1.3.0.6
	OPENBSD_5_7_BASE:1.3
	OPENBSD_5_6:1.3.0.4
	OPENBSD_5_6_BASE:1.3
	OPENBSD_5_5:1.3.0.2
	OPENBSD_5_5_BASE:1.3;
locks; strict;
comment	@ * @;


1.9
date	2017.05.27.13.53.07;	author ajacoutot;	state Exp;
branches;
next	1.8;
commitid	SDpZdO8hfYIlMQ9j;

1.8
date	2015.12.21.10.17.32;	author ajacoutot;	state Exp;
branches;
next	1.7;
commitid	mtNXgqH8qNXf2w4k;

1.7
date	2015.12.21.09.32.01;	author ajacoutot;	state Exp;
branches;
next	1.6;
commitid	89yBmxf8zl7f4edS;

1.6
date	2015.12.20.09.38.14;	author ajacoutot;	state Exp;
branches;
next	1.5;
commitid	SUL8g4c0gjEg4ZTf;

1.5
date	2015.12.19.14.13.45;	author ajacoutot;	state Exp;
branches;
next	1.4;
commitid	jNmUMRLbxGUrj246;

1.4
date	2015.04.13.13.30.02;	author zhuk;	state Exp;
branches;
next	1.3;
commitid	zsGVZxArmIZNgFcF;

1.3
date	2014.01.30.16.58.58;	author ajacoutot;	state Exp;
branches;
next	1.2;

1.2
date	2014.01.29.16.49.33;	author zhuk;	state Exp;
branches;
next	1.1;

1.1
date	2014.01.28.20.50.40;	author zhuk;	state Exp;
branches;
next	;


desc
@@


1.9
log
@Cope with recent rc.subr(8) change (removal of rc_bg etc.)
@
text
@#!/bin/sh
#
# $OpenBSD: kdm.rc,v 1.8 2015/12/21 10:17:32 ajacoutot Exp $

daemon="${TRUEPREFIX}/bin/kdm"

. /etc/rc.d/rc.subr

rc_reload=NO

if [ -n "${INRC}" ]; then
# on boot: wait for ttys to be initialized in bg so we we don't hang in _rc_wait
rc_start() {
	( local i=0
	while ! pgrep -qf "^/usr/libexec/getty "; do
		sleep 1
		[ $((i++)) -ge 10 ] && return 1
	done
	# XXX KDM fails to receive input unless getty is given
	# some more time for initialization.
	sleep 2
	${rcexec} "${daemon} ${daemon_flags}" ) &
	rc_check() { return 0; } # ignore status, we're waiting for getty
}
fi

rc_cmd $1
@


1.8
log
@Override _rc_wait(); these scripts are a bit hackish because they need
/usr/libexec/getty to be started but it is started later in the boot
process.
@
text
@d3 1
a3 1
# $OpenBSD: kdm.rc,v 1.7 2015/12/21 09:32:01 ajacoutot Exp $
d12 1
a12 5
# on boot: make sure we don't hang in _rc_wait
_rc_wait() {
	return 0
}
# on boot: wait for ttys to be initialized
d23 1
@


1.7
log
@Revert, this breaks rc_stop() for these daemons.
I'm working on a better fix.

spotted by dcoppa@@
@
text
@d3 1
a3 1
# $OpenBSD: kdm.rc,v 1.6 2015/12/20 09:38:14 ajacoutot Exp $
a10 2
# on boot: wait for ttys to be initialized and make sure we don't hang in
# "_rc_wait start" by hacking pexp
d12 5
@


1.6
log
@Move pexp where it belongs.
@
text
@d3 1
a3 1
# $OpenBSD: kdm.rc,v 1.5 2015/12/19 14:13:45 ajacoutot Exp $
a14 1
	pexp="/sbin/init"
@


1.5
log
@Cope with new rc.subr behavior so that we don't end hanging in _rc_wait
start.
@
text
@d3 1
a3 1
# $OpenBSD: kdm.rc,v 1.4 2015/04/13 13:30:02 zhuk Exp $
d14 1
a15 1
rc_start() {
@


1.4
log
@Add an XXX comment about "sleep 2" in rc script, requested by aja@@.
@
text
@d3 1
a3 1
# $OpenBSD: kdm.rc,v 1.3 2014/01/30 16:58:58 ajacoutot Exp $
d11 2
a12 1
# on boot: wait for ttys to be initialized
d14 1
@


1.3
log
@Sync rc_start() between these 3 daemons since they do the same.
Keeping the sleep 2 in kdm.rc -- not sure why this is needed, zhuk@@ probably
knows.
@
text
@d3 1
a3 1
# $OpenBSD: kdm.rc,v 1.2 2014/01/29 16:49:33 zhuk Exp $
d19 2
@


1.2
log
@Simplify the rc script by conditional define of function.

Idea and initial implementation (in gdm) from ajacoutot@@.
@
text
@d3 1
a3 1
# $OpenBSD: kdm.rc,v 1.1 2014/01/28 20:50:40 zhuk Exp $
d11 1
@


1.1
log
@Polish KDM support:

  * Add rc script. Just add "kdm" to pkg_scripts in rc.conf.local now.

  * Add a dirty hack for console-kit-daemon until ConsoleKit gets fixed.
    We don't support multi-seat X now anyway.

Please note: while the rc script is based on gdm.rc, it has somewhat changed
logic. So if you're using KDM, please test it in rc.conf.local before we
get to 5.5 release, to make sure it doesn't break. It shouldn't, of course,
but, you know, IT happens.
@
text
@d3 1
a3 1
# $OpenBSD$
d11 1
d13 7
a19 14
	if pgrep -qf "^/usr/libexec/getty "; then
		${rcexec} "${daemon} ${daemon_flags}"
	else
		# OS is starting, workaround like in gdm.rc
		( local i=0
		while [[ $i -lt 10 ]]; do
			pgrep -qf "^/usr/libexec/getty " && break
			sleep 1
			i=$((i+1))
		done
		[[ $i -ge 10 ]] && return 1
		sleep 2
		${rcexec} "${daemon} ${daemon_flags}" ) &
	fi
d21 1
@

