# cvs/ports/net/isc-bind/pkg/isc_named.rc,v
head	1.3;
access;
symbols
	OPENBSD_6_0:1.3.0.8
	OPENBSD_6_0_BASE:1.3
	OPENBSD_5_9:1.3.0.4
	OPENBSD_5_9_BASE:1.3
	OPENBSD_5_8:1.3.0.6
	OPENBSD_5_8_BASE:1.3
	OPENBSD_5_7:1.3.0.2
	OPENBSD_5_7_BASE:1.3
	OPENBSD_5_6:1.2.0.6
	OPENBSD_5_6_BASE:1.2
	OPENBSD_5_5:1.2.0.4
	OPENBSD_5_5_BASE:1.2
	OPENBSD_5_4:1.2.0.2
	OPENBSD_5_4_BASE:1.2
	OPENBSD_5_3:1.1.1.1.0.4
	OPENBSD_5_3_BASE:1.1.1.1
	OPENBSD_5_2:1.1.1.1.0.2
	OPENBSD_5_2_BASE:1.1.1.1
	sthen_20120719:1.1.1.1
	sthen:1.1.1;
locks; strict;
comment	@ * @;


1.3
date	2014.08.22.20.32.13;	author sthen;	state Exp;
branches;
next	1.2;
commitid	yw5iSs0sfRNvewuY;

1.2
date	2013.04.03.20.17.23;	author sthen;	state Exp;
branches;
next	1.1;

1.1
date	2012.07.19.08.36.45;	author sthen;	state Exp;
branches
	1.1.1.1;
next	;

1.1.1.1
date	2012.07.19.08.36.45;	author sthen;	state Exp;
branches
	1.1.1.1.4.1;
next	;

1.1.1.1.4.1
date	2013.05.01.09.06.03;	author sthen;	state Exp;
branches;
next	;


desc
@@


1.3
log
@install sample configs and switch to running as _bind uid
@
text
@#!/bin/sh
#
# $OpenBSD: isc_named.rc,v 1.2 2013/04/03 20:17:23 sthen Exp $

daemon="${TRUEPREFIX}/sbin/named"
daemon_flags="-t ${LOCALSTATEDIR}/named -u _bind -U 4"

. /etc/rc.d/rc.subr

rc_pre() {
	if ! cmp -s ${SYSCONFDIR}/rndc.key ${LOCALSTATEDIR}/named/etc/rndc.key ; then
		if ${TRUEPREFIX}/sbin/rndc-confgen -a -t ${LOCALSTATEDIR}/named \
		    >/dev/null 2>&1; then
			chmod 0640 ${LOCALSTATEDIR}/named/etc/rndc.key \
			    >/dev/null 2>&1
		fi
	fi
}

rc_cmd $1
@


1.2
log
@Normally newer versions of BIND default to one UDP listener per interface
per CPU. As found by Marc Peters, this doesn't work too well on a t5120
with 64 threads, so change the default settings in the rc.d script to -U 4
to cap this to 4, or the number of CPUs if less.

As usual with rc.d scripts, if you need to override flags, set
isc_named_flags="..." in rc.conf.local.
@
text
@d3 1
a3 1
# $OpenBSD: isc_named.rc,v 1.1.1.1 2012/07/19 08:36:45 sthen Exp $
d6 1
a6 1
daemon_flags="-t ${LOCALSTATEDIR}/named -u named -U 4"
d11 1
a11 1
	if ! cmp -s /etc/rndc.key ${LOCALSTATEDIR}/named/etc/rndc.key ; then
@


1.1
log
@Initial revision
@
text
@d3 1
a3 1
# $OpenBSD$
d6 1
a6 1
daemon_flags="-t ${LOCALSTATEDIR}/named -u named"
@


1.1.1.1
log
@import BIND 9.9.1-P1 as ports/net/isc-bind - this is newer than the
version of BIND than is in the base OS (some people require features
from this version e.g. DNS64), but note that it does not include
the hardening changes made to the version in base.

feedback from naddy@@ giovanni@@, ok giovanni@@.

"BIND is open source software that implements the Domain Name System
(DNS) protocols for the Internet. It is a reference implementation
of those protocols, but it is also production-grade software,
suitable for use in high-volume and high-reliability applications."
@
text
@@


1.1.1.1.4.1
log
@MFC - BIND 9.9.2-P2: CVE-2013-2266 https://kb.isc.org/article/AA-00871
@
text
@d3 1
a3 1
# $OpenBSD: isc_named.rc,v 1.2 2013/04/03 20:17:23 sthen Exp $
d6 1
a6 1
daemon_flags="-t ${LOCALSTATEDIR}/named -u named -U 4"
@

