# cvs/ports/sysutils/mcollective/pkg/mcollectived.rc,v
head	1.3;
access;
symbols
	OPENBSD_6_0:1.2.0.14
	OPENBSD_6_0_BASE:1.2
	OPENBSD_5_9:1.2.0.10
	OPENBSD_5_9_BASE:1.2
	OPENBSD_5_8:1.2.0.12
	OPENBSD_5_8_BASE:1.2
	OPENBSD_5_7:1.2.0.8
	OPENBSD_5_7_BASE:1.2
	OPENBSD_5_6:1.2.0.6
	OPENBSD_5_6_BASE:1.2
	OPENBSD_5_5:1.2.0.4
	OPENBSD_5_5_BASE:1.2
	OPENBSD_5_4:1.2.0.2
	OPENBSD_5_4_BASE:1.2
	OPENBSD_5_3:1.1.0.6
	OPENBSD_5_3_BASE:1.1
	OPENBSD_5_2:1.1.0.4
	OPENBSD_5_2_BASE:1.1
	OPENBSD_5_1_BASE:1.1
	OPENBSD_5_1:1.1.0.2;
locks; strict;
comment	@ * @;


1.3
date	2016.09.21.20.27.40;	author sebastia;	state Exp;
branches;
next	1.2;
commitid	9fKZMhKJ3REFaIIE;

1.2
date	2013.03.13.22.54.26;	author ajacoutot;	state Exp;
branches;
next	1.1;

1.1
date	2011.09.20.14.03.20;	author uwe;	state Exp;
branches;
next	;


desc
@@


1.3
log
@mcollectived needs --pid parameter in order to properly daemonize

feedback/OK jasper@@, ajacoutot@@
@
text
@#!/bin/sh
#
# $OpenBSD: mcollectived.rc,v 1.2 2013/03/13 22:54:26 ajacoutot Exp $

daemon="${TRUEPREFIX}/sbin/mcollectived"
daemon_flags="--pid=/var/run/mcollectived.pid"

. /etc/rc.d/rc.subr

pexp="${RUBY} ${daemon}${daemon_flags:+ ${daemon_flags}}"
rc_reload=NO

rc_cmd $1
@


1.2
log
@Use bundled modules since they are certified to work with that particular
version of mcollective.
Force ruby 1.8; version 1.9 is not functionnal yet.
Simplify patching the config files by using a regex.
Use BSD_INSTALL_*
Use $RUBY in the rc script.
Take maintainer.

More work to come, so still not hooked to the build.
@
text
@d3 1
a3 1
# $OpenBSD: mcollectived.rc,v 1.1 2011/09/20 14:03:20 uwe Exp $
d6 1
@


1.1
log
@Import ruby-mcollective-1.3.0

This port is based on the port source that is currently bundled with
MCollective but which has not been updated for ages.

Homepage: http://puppetlabs.com/mcollective/
Source: https://github.com/puppetlabs/marionette-collective
@
text
@d3 1
a3 1
# $OpenBSD$
d9 1
a9 1
pexp=".*ruby.* ${daemon}${daemon_flags:+ ${daemon_flags}}"
@

