# cvs/ports/net/quagga/pkg/quagga_pimd.rc,v
head	1.1;
access;
symbols
	OPENBSD_6_0:1.1.0.2
	OPENBSD_6_0_BASE:1.1;
locks; strict;
comment	@ * @;


1.1
date	2016.03.22.12.57.54;	author jca;	state Exp;
branches;
next	;
commitid	UMHjWDfFsEz4krLO;


desc
@@


1.1
log
@Update to 1.0.20160315, ok sthen@@ (maintainer)

Announce:

  https://lists.quagga.net/pipermail/quagga-users/2016-March/014293.html

ChangeLog:

  http://savannah.spinellicreations.com/quagga/quagga-1.0.20160315.changelog.txt
  http://savannah.spinellicreations.com/quagga/quagga-1.0.20160309.changelog.txt

Note that this update brings a fix for CVE-2016-2342.  This vulnerability
affects quagga MPLS VPN setups which probably aren't possible on OpenBSD.
@
text
@#!/bin/sh
#
# $OpenBSD$

daemon="${TRUEPREFIX}/sbin/pimd -d"

. /etc/rc.d/rc.subr

rc_cmd $1
@
