# cvs/ports/net/ptpd/pkg/ptpd.rc,v
head	1.1;
access;
symbols
	OPENBSD_6_2:1.1.0.2
	OPENBSD_6_2_BASE:1.1;
locks; strict;
comment	@ * @;


1.1
date	2017.09.09.21.22.49;	author naddy;	state Exp;
branches;
next	;
commitid	JfjQznIXnlRBDFDh;


desc
@@


1.1
log
@* Update to 1.1.0.
* Replace MESSAGE with rc file.
* Drop maintainership.
* Set PORTROACH.

ptpd-1.x implements PTPv1 (IEEE 1588-2002)
ptpd-2.x implements PTPv2 (IEEE 1588-2008)
PTPv2 is not backward compatible with PTPv1.
@
text
@#!/bin/sh
#
# $OpenBSD$

daemon="${TRUEPREFIX}/sbin/ptpd"

. /etc/rc.d/rc.subr

rc_cmd $1
@
