# cvs/ports/mail/rmilter/pkg/rmilter.rc,v
head	1.1;
branch	1.1.1;
access;
symbols
	OPENBSD_6_2:1.1.1.1.0.4
	OPENBSD_6_2_BASE:1.1.1.1
	OPENBSD_6_1:1.1.1.1.0.2
	OPENBSD_6_1_BASE:1.1.1.1
	sthen_20161202:1.1.1.1
	sthen:1.1.1;
locks; strict;
comment	@ * @;


1.1
date	2016.12.02.14.57.41;	author sthen;	state Exp;
branches
	1.1.1.1;
next	;
commitid	dd8meP197Xg3dwF8;

1.1.1.1
date	2016.12.02.14.57.41;	author sthen;	state Exp;
branches;
next	;
commitid	dd8meP197Xg3dwF8;


desc
@@



1.1
log
@Initial revision
@
text
@#!/bin/sh
#
# $OpenBSD$

daemon="${TRUEPREFIX}/sbin/rmilter"
daemon_flags="-c ${SYSCONFDIR}/rmilter/rmilter.conf"
daemon_user="_rmilter"

. /etc/rc.d/rc.subr

rc_pre() {
	/usr/bin/install -d -o ${daemon_user} /var/run/rmilter
}

rc_reload() {
        pkill -USR1 -T "${daemon_rtable}" -xf "${pexp}"
}

rc_cmd $1
@


1.1.1.1
log
@import ports/mail/rmilter, testing/tweaks Johan Huldtgren, review/ok jca

Rmilter is used to integrate rspamd with any milter-compatible MTA, for
example postfix or sendmail. Features include:

- ClamAV scanning (via unix or tcp socket).
- Rspamd scanning
- Greylisting and rate-limiting with redis upstream
- Auto-whitelisting (internal and via redis upstream)
- Replies check (whitelisting replies to sent messages)
- Passing messages and/or their headers to beanstalk servers
@
text
@@
