# cvs/ports/www/nginx/pkg/Attic/PLIST,v
head	1.17;
access;
symbols
	OPENBSD_6_0:1.16.0.8
	OPENBSD_6_0_BASE:1.16
	OPENBSD_5_9:1.16.0.4
	OPENBSD_5_9_BASE:1.16
	OPENBSD_5_8:1.16.0.6
	OPENBSD_5_8_BASE:1.16
	OPENBSD_5_7:1.16.0.2
	OPENBSD_5_7_BASE:1.16
	OPENBSD_5_2:1.9.0.2
	OPENBSD_5_2_BASE:1.9
	OPENBSD_5_1_BASE:1.8
	OPENBSD_5_1:1.8.0.4
	OPENBSD_5_0:1.8.0.2
	OPENBSD_5_0_BASE:1.8
	OPENBSD_4_9:1.7.0.2
	OPENBSD_4_9_BASE:1.7
	OPENBSD_4_8:1.4.0.8
	OPENBSD_4_8_BASE:1.4
	OPENBSD_4_7:1.4.0.6
	OPENBSD_4_7_BASE:1.4
	OPENBSD_4_6:1.4.0.4
	OPENBSD_4_6_BASE:1.4
	OPENBSD_4_5:1.4.0.2
	OPENBSD_4_5_BASE:1.4
	OPENBSD_4_4:1.3.0.2
	OPENBSD_4_4_BASE:1.3
	OPENBSD_4_3:1.2.0.2
	OPENBSD_4_3_BASE:1.2
	merdely_20070904:1.1.1.1
	merdely:1.1.1;
locks; strict;
comment	@# @;


1.17
date	2016.11.04.22.13.43;	author jeremy;	state dead;
branches;
next	1.16;
commitid	Oee7bZXxLwPn3r9t;

1.16
date	2014.10.23.12.13.12;	author sthen;	state Exp;
branches;
next	1.15;
commitid	jQCiePNARms9IfsU;

1.15
date	2014.09.09.14.05.48;	author sthen;	state Exp;
branches;
next	1.14;
commitid	w2JJgh9xh2JAWyXg;

1.14
date	2014.09.09.13.48.28;	author sthen;	state Exp;
branches;
next	1.13;
commitid	vtDVj30kfNFKowNU;

1.13
date	2014.08.27.06.46.49;	author ajacoutot;	state Exp;
branches;
next	1.12;
commitid	EgWPCQ6ECli0Ppnk;

1.12
date	2014.08.26.19.48.06;	author robert;	state Exp;
branches;
next	1.11;
commitid	Rxd8gSExVfEEQmSV;

1.11
date	2014.08.26.19.43.59;	author robert;	state Exp;
branches;
next	1.10;
commitid	Y2vy3KpSALTztMGt;

1.10
date	2012.09.11.16.01.45;	author sthen;	state dead;
branches;
next	1.9;

1.9
date	2012.02.19.15.02.08;	author robert;	state Exp;
branches;
next	1.8;

1.8
date	2011.05.17.03.31.00;	author william;	state Exp;
branches;
next	1.7;

1.7
date	2010.11.10.08.18.18;	author sthen;	state Exp;
branches;
next	1.6;

1.6
date	2010.11.04.07.39.44;	author robert;	state Exp;
branches;
next	1.5;

1.5
date	2010.11.04.01.41.11;	author dlg;	state Exp;
branches;
next	1.4;

1.4
date	2008.09.01.15.47.44;	author sthen;	state Exp;
branches;
next	1.3;

1.3
date	2008.03.29.13.50.39;	author bernd;	state Exp;
branches;
next	1.2;

1.2
date	2007.12.31.02.20.01;	author okan;	state Exp;
branches;
next	1.1;

1.1
date	2007.09.04.21.58.26;	author merdely;	state Exp;
branches
	1.1.1.1;
next	;

1.1.1.1
date	2007.09.04.21.58.26;	author merdely;	state Exp;
branches;
next	;


desc
@@


1.17
log
@Update nginx to use new dynamic modules feature

Switch all extrnal modules to be dynamic modules, and update them to use the
latest versions.  Use separate subpackages for all dynamic modules, so you
can easily choose which nginx features you want without rebuilding nginx
with a specific combination of FLAVORs.

Switch all internal modules that can be built dynamically to use dynamic
modules.

Start building the image_filter, geoip, perl, and xlst modules as dynamic
modules.

This requires changes to your nginx configuration if you were previously
using a FLAVORed version of nginx, or using the mail or stream features.

OK robert@@ (MAINTAINER)
positive feedback from william@@ (MAINTAINER), landry@@, sthen@@
@
text
@@@comment $OpenBSD: PLIST,v 1.16 2014/10/23 12:13:12 sthen Exp $
@@pkgpath www/nginx/stable${FLAVOR_COMMA}
@@pkgpath www/nginx/devel${FLAVOR_COMMA}
@@man man/man8/nginx.8
@@bin sbin/nginx
@@mode 700
@@owner www
@@group www
@@sample ${NGINX_DIR}/tmp/
@@mode
@@owner
@@group
share/doc/pkg-readmes/${FULLPKGNAME}
share/nginx/
@@sample ${SYSCONFDIR}/nginx/
share/nginx/fastcgi_params
@@sample ${SYSCONFDIR}/nginx/fastcgi_params
share/nginx/html/
@@sample ${NGINX_DIR}/html/
share/nginx/html/50x.html
@@sample ${NGINX_DIR}/html/50x.html
share/nginx/html/index.html
@@sample ${NGINX_DIR}/html/index.html
share/nginx/koi-utf
@@sample ${SYSCONFDIR}/nginx/koi-utf
share/nginx/koi-win
@@sample ${SYSCONFDIR}/nginx/koi-win
share/nginx/mime.types
@@sample ${SYSCONFDIR}/nginx/mime.types
share/nginx/nginx.conf
@@sample ${SYSCONFDIR}/nginx/nginx.conf
share/nginx/scgi_params
@@sample ${SYSCONFDIR}/nginx/scgi_params
share/nginx/uwsgi_params
@@sample ${SYSCONFDIR}/nginx/uwsgi_params
share/nginx/win-utf
@@sample ${SYSCONFDIR}/nginx/win-utf
@@rcscript ${RCDIR}/nginx
@


1.16
log
@Use specific flavour extensions for @@pkgpath lines in each package, so that
pkg_add -u is handled without prompting each time, as noticed by myself and
jasper. Modelled on similar in mail/mutt.  ok jasper@@ espie@@
@
text
@d1 1
a1 1
@@comment $OpenBSD: PLIST,v 1.15 2014/09/09 14:05:48 sthen Exp $
@


1.15
log
@@@sample nginx config files, suggested by/ok aja@@
@
text
@d1 3
a3 5
@@comment $OpenBSD: PLIST,v 1.14 2014/09/09 13:48:28 sthen Exp $
@@pkgpath www/nginx/stable
@@pkgpath www/nginx/devel
%%passenger%%
!%%passenger%%
@


1.14
log
@add a README to www/nginx with clues about log-rotation, ok aja@@
@
text
@d1 1
a1 1
@@comment $OpenBSD: PLIST,v 1.13 2014/08/27 06:46:49 ajacoutot Exp $
d17 1
d19 1
d27 1
d29 1
d31 1
d33 1
d35 1
d37 1
d39 1
@


1.13
log
@Drop out-dated README.

ok robert@@
@
text
@d1 1
a1 1
@@comment $OpenBSD: PLIST,v 1.12 2014/08/26 19:48:06 robert Exp $
d15 1
@


1.12
log
@add @@pkgpaths
@
text
@d1 1
a1 1
@@comment $OpenBSD: PLIST,v 1.11 2014/08/26 19:43:59 robert Exp $
a14 1
share/doc/pkg-readmes/${FULLPKGNAME}
@


1.11
log
@- remove the devel nginx packages
- apply chroot and other patches from the base system
- update to 1.7.4
@
text
@d1 3
a3 1
@@comment $OpenBSD: PLIST,v 1.1 2012/09/11 16:01:45 sthen Exp $
@


1.10
log
@Update nginx stable in ports to 1.2.x and add devel (1.3.x), this is
separate from the version in the base OS which deliberately does not include
the mail proxy and additional modules. rc.d/enginx script is added (not
rc.d/nginx to avoid conflicting with base). Lua module is also supported.
Most of the work done by william@@, with a few tweaks by me and ok william.

README could use more work and as pointed out by ajacoutot we need to
review permissions of the tmp directory, but the diff is already long enough.
@@pkgpath markers are set to update from the 1.0.x packages to 1.2.x (with or
without the passenger flavour as appropriate).
@
text
@d1 3
a3 1
@@comment $OpenBSD: PLIST,v 1.9 2012/02/19 15:02:08 robert Exp $
d13 1
d29 1
@


1.9
log
@do not handle the same files that are handled by the in-tree nginx
and make the port use the www user so that both nginx can be used
on the same system if needed
@
text
@d1 1
a1 1
@@comment $OpenBSD: PLIST,v 1.8 2011/05/17 03:31:00 william Exp $
@


1.8
log
@Update to nginx-1.0.2

Originally from Piotr for 1.0.0, thanks

ok jasper
@
text
@d1 1
a1 3
@@comment $OpenBSD: PLIST,v 1.7 2010/11/10 08:18:18 sthen Exp $
@@newgroup _nginx:593
@@newuser _nginx:593:_nginx:daemon:nginx user:/nonexistent:/sbin/nologin
a3 3
@@mode 755
@@sample ${NGINX_DIR}/
@@sample /var/log/nginx/
d5 2
a6 2
@@owner _nginx
@@group _nginx
a10 1
@@sample ${SYSCONFDIR}/nginx/
a12 1
@@sample ${SYSCONFDIR}/nginx/fastcgi_params
a19 1
@@sample ${SYSCONFDIR}/nginx/koi-utf
a20 1
@@sample ${SYSCONFDIR}/nginx/koi-win
a21 1
@@sample ${SYSCONFDIR}/nginx/mime.types
a22 1
@@sample ${SYSCONFDIR}/nginx/nginx.conf
a23 1
@@sample ${SYSCONFDIR}/nginx/scgi_params
a24 1
@@sample ${SYSCONFDIR}/nginx/uwsgi_params
a25 2
@@sample ${SYSCONFDIR}/nginx/win-utf
@@rcscript ${RCDIR}/nginx
@


1.7
log
@- Install files for the SCGI and uwsgi modules, from Piotr Sikora.
makes sense to dlg@@, ok william@@ (maintainer).

- Fix bad REVISION bump in an .if ${FLAVOR... block from last commit.
@
text
@d1 1
a1 1
@@comment $OpenBSD: PLIST,v 1.6 2010/11/04 07:39:44 robert Exp $
d4 1
@


1.6
log
@Add an rc script for nginx too
@
text
@d1 1
a1 1
@@comment $OpenBSD: PLIST,v 1.5 2010/11/04 01:41:11 dlg Exp $
d33 4
@


1.5
log
@update nginx to 0.8.53. the main flavor builds and works fine. the
passenger build works, but i havent got a test case for it.

ok william@@
@
text
@d1 1
a1 1
@@comment $OpenBSD: PLIST,v 1.4 2008/09/01 15:47:44 sthen Exp $
d16 1
d35 1
@


1.4
log
@Maintenance update to latest stable version 0.6.32, from maintainer
Darrin Chandler on ports@@.

"Thanks to Toni Mueller for prodding, and to Simon Kuhnle for prodding,
testing, and the patch (I had the same diff just before lock but kept
forgetting after unlock)." -dwc
@
text
@d1 1
a1 1
@@comment $OpenBSD: PLIST,v 1.3 2008/03/29 13:50:39 bernd Exp $
a14 1
share/nginx/
@


1.3
log
@Give nginx a chance to write stuff into /var/nginx/tmp.

ok Darrin Chandler (MAINTAINER)
@
text
@d1 1
a1 1
@@comment $OpenBSD: PLIST,v 1.2 2007/12/31 02:20:01 okan Exp $
d4 1
a4 1
sbin/nginx
@


1.2
log
@update from MAINTAINER; plus a bunch of tweaks from me

ok merdely@@
@
text
@d1 1
a1 1
@@comment $OpenBSD: PLIST,v 1.1.1.1 2007/09/04 21:58:26 merdely Exp $
d9 2
d13 2
@


1.1
log
@Initial revision
@
text
@d1 1
a1 1
@@comment $OpenBSD$
d3 1
a3 1
@@newuser _nginx:593:_nginx:daemon:nginx user:/nonexistant:/sbin/nologin
@


1.1.1.1
log
@Import nginx 0.5.31

nginx [engine x] is a HTTP server and IMAP/POP3 proxy server.

The basic HTTP features:
  * Handling of the static files, index files, and autoindexing
  * Accelerated reverse proxying without caching, simple load 
    balancing and fault tolerance
  * Accelerated support without caching of the remote FastCGI
    servers, simple load balancing and fault tolerance
  * Modular architecture, filters including gzipping, byte
    ranges, chunked responses, and SSI-filter; Several subrequests
    in one page handling in SSI-filter via FastCGI or proxy
    running in parallel
  * the SSL support

The IMAP/POP3 proxy server features:
  * User redirection to IMAP/POP3 backend using an external HTTP
    authentication server
  * Plain text authentication (LOGIN, USER/PASS)
  * SSL and STARTTLS support

ok and help steven@@
@
text
@@
