# cvs/ports/sysutils/syslog-ng/pkg/syslog_ng.rc,v
head	1.6;
access;
symbols
	OPENBSD_6_1:1.6.0.14
	OPENBSD_6_1_BASE:1.6
	OPENBSD_6_0:1.6.0.12
	OPENBSD_6_0_BASE:1.6
	OPENBSD_5_9:1.6.0.8
	OPENBSD_5_9_BASE:1.6
	OPENBSD_5_8:1.6.0.10
	OPENBSD_5_8_BASE:1.6
	OPENBSD_5_7:1.6.0.6
	OPENBSD_5_7_BASE:1.6
	OPENBSD_5_6:1.6.0.4
	OPENBSD_5_6_BASE:1.6
	OPENBSD_5_5:1.6.0.2
	OPENBSD_5_5_BASE:1.6
	OPENBSD_5_4:1.5.0.4
	OPENBSD_5_4_BASE:1.5
	OPENBSD_5_3:1.5.0.2
	OPENBSD_5_3_BASE:1.5
	OPENBSD_5_2:1.4.0.2
	OPENBSD_5_2_BASE:1.4
	OPENBSD_5_1_BASE:1.1
	OPENBSD_5_1:1.1.0.2;
locks; strict;
comment	@ * @;


1.6
date	2013.11.27.07.52.54;	author ajacoutot;	state Exp;
branches;
next	1.5;

1.5
date	2012.09.19.12.25.40;	author ajacoutot;	state Exp;
branches;
next	1.4;

1.4
date	2012.06.08.21.06.46;	author steven;	state Exp;
branches;
next	1.3;

1.3
date	2012.05.28.16.31.31;	author ajacoutot;	state Exp;
branches;
next	1.2;

1.2
date	2012.05.23.06.44.50;	author ajacoutot;	state Exp;
branches;
next	1.1;

1.1
date	2011.09.14.02.34.06;	author william;	state Exp;
branches;
next	;


desc
@@


1.6
log
@Drop the INRC special case -- switching syslogd for syslog-ng should be
left as an exercice for whoever needs/wants to do it.

ok sthen@@ steven@@
@
text
@#!/bin/sh
#
# $OpenBSD: syslog_ng.rc,v 1.5 2012/09/19 12:25:40 ajacoutot Exp $

daemon="${TRUEPREFIX}/sbin/syslog-ng"
daemon_flags="-p /var/run/syslog.pid"

. /etc/rc.d/rc.subr

pexp="${daemon}"

rc_cmd $1
@


1.5
log
@Unfuck ld.so.hints when run late in the boot process.

prodded by and ok espie@@
@
text
@d3 1
a3 1
# $OpenBSD: syslog_ng.rc,v 1.4 2012/06/08 21:06:46 steven Exp $
a10 6

# needed in case syslog-ng(8) is started early in the boot sequence
# before rc(8) had a chance to run ldconfig
rc_pre() {
	[ "${INRC}" ] && /sbin/ldconfig -m ${LOCALBASE}/lib || return 0
}
@


1.4
log
@s/PREFIX/TRUEPREFIX/ for consistency.

from robert at peichaer.org

ok ajacoutot@@
@
text
@d3 1
a3 1
# $OpenBSD: syslog_ng.rc,v 1.2 2012/05/23 06:44:50 ajacoutot Exp $
d15 1
a15 1
	[ "${INRC}" ] && /sbin/ldconfig ${LOCALBASE}/lib || return 0
@


1.3
log
@Set the PID to /var/run/syslog.pid by default to make syslog-ng happy
with newsyslog(8); this can of course by modified using syslog_ng_flags.
While here, run syslog-ng.conf(5) through tbl(1) (from sthen@@).

discussed with and ok sthen@@ steven@@ (maintainer)
@
text
@d5 1
a5 1
daemon="${PREFIX}/sbin/syslog-ng"
@


1.2
log
@If syslog-ng is started early during boot-up (that is if one replaces
the syslogd rc.d(8) script with the syslog_ng one), it will fail to
launch because ldconfig(8) hasn't had a chance to run yet. Add an
rc_pre() that will take care of that.

ok sthen@@ steven@@ (maintainer)
@
text
@d3 1
a3 1
# $OpenBSD: syslog_ng.rc,v 1.1 2011/09/14 02:34:06 william Exp $
d6 1
d9 2
@


1.1
log
@add rcscript

ok ajacoutot steven (MAINTAINER)
@
text
@d3 1
a3 1
# $OpenBSD$
d8 6
@

