# cvs/ports/mail/mimedefang/pkg/Attic/mimedefang_daemon.rc,v
head	1.8;
access;
symbols
	OPENBSD_5_2:1.6.0.2
	OPENBSD_5_2_BASE:1.6
	OPENBSD_5_1_BASE:1.5
	OPENBSD_5_1:1.5.0.6
	OPENBSD_5_0:1.5.0.4
	OPENBSD_5_0_BASE:1.5
	OPENBSD_4_9:1.5.0.2
	OPENBSD_4_9_BASE:1.5;
locks; strict;
comment	@ * @;


1.8
date	2012.12.05.18.51.56;	author ajacoutot;	state dead;
branches;
next	1.7;

1.7
date	2012.09.18.15.34.39;	author ajacoutot;	state Exp;
branches;
next	1.6;

1.6
date	2012.05.31.18.31.42;	author ajacoutot;	state Exp;
branches;
next	1.5;

1.5
date	2010.12.27.14.50.23;	author ajacoutot;	state Exp;
branches;
next	1.4;

1.4
date	2010.12.24.10.40.05;	author ajacoutot;	state Exp;
branches;
next	1.3;

1.3
date	2010.12.11.12.57.37;	author ajacoutot;	state Exp;
branches;
next	1.2;

1.2
date	2010.12.09.01.44.15;	author okan;	state Exp;
branches;
next	1.1;

1.1
date	2010.12.07.02.56.19;	author okan;	state Exp;
branches;
next	;


desc
@@


1.8
log
@Add `-r' to mimedefang flags by default.
Add an example for filter_relay() that works out of the box.
Rename the mimedefang_daemon rc.d(8) script to mimedefang as it's
supposed to be...

ok okan@@ (maintainer)
@
text
@#!/bin/sh
#
# $OpenBSD: mimedefang_daemon.rc,v 1.7 2012/09/18 15:34:39 ajacoutot Exp $

daemon="${TRUEPREFIX}/bin/mimedefang"
daemon_flags="-G -U _mdefang -P /var/spool/MIMEDefang/mimedefang.pid -p /var/spool/MIMEDefang/mimedefang.sock -m /var/spool/MIMEDefang/mimedefang-multiplexor.sock"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
@


1.7
log
@* set default email to postmaster@@localhost since this alias exists by default
and change AdminName accordingly
* disable the filter_bad_filename filter by default, it is the admin's
job to choose whether he wants to enable it or not
* remove useless spamassassin.cf patch
* run mimedefang with -G so that external helpers (e.g. antivirus scanners)
can access its socket by being added to the _mdefang group

ok okan@@ (maintainer)
@
text
@d3 1
a3 1
# $OpenBSD: mimedefang_daemon.rc,v 1.6 2012/05/31 18:31:42 ajacoutot Exp $
@


1.6
log
@Update to mimedefang-2.73.

Rework the port to comply with upstream defaults and generally simplify
it. Set default options in rc.d(8) following the default init script
behavior.
Activate ipheader and install gen-ip-validator.pl.

ok okan@@ (maintainer)
@
text
@d3 1
a3 1
# $OpenBSD: mimedefang_daemon.rc,v 1.5 2010/12/27 14:50:23 ajacoutot Exp $
d6 1
a6 1
daemon_flags="-U _mdefang -P /var/spool/MIMEDefang/mimedefang.pid -p /var/spool/MIMEDefang/mimedefang.sock -m /var/spool/MIMEDefang/mimedefang-multiplexor.sock"
@


1.5
log
@Simplify after recent rc.subr change.
The framework is now stable and we will start documenting it (at last).
@
text
@d3 1
a3 1
# $OpenBSD: mimedefang_daemon.rc,v 1.4 2010/12/24 10:40:05 ajacoutot Exp $
d6 1
a6 1
daemon_flags="-U _mdefang -p /var/run/mdefang/mimedefang.socket -m /var/run/mdefang/multiplexor.socket"
a10 4

rc_pre() {
	install -d -o _mdefang -g _mdefang /var/run/mdefang
}
@


1.4
log
@Cope with recent rc.subr changes.
@
text
@d3 1
a3 3
# $OpenBSD: mimedefang_daemon.rc,v 1.3 2010/12/11 12:57:37 ajacoutot Exp $

. /etc/rc.d/rc.subr
d8 1
a8 1
rc_conf
@


1.3
log
@Move these to rc_reload=NO.
@
text
@d3 1
a3 1
# $OpenBSD: mimedefang_daemon.rc,v 1.2 2010/12/09 01:44:15 okan Exp $
d9 3
@


1.2
log
@add an appropriate rc_reload()

proded by and ok ajacoutot@@
@
text
@d3 1
a3 1
# $OpenBSD: mimedefang_daemon.rc,v 1.1 2010/12/07 02:56:19 okan Exp $
d9 1
a12 4
}

rc_reload() {
	rc_err "$0: reload is not supported"
@


1.1
log
@- update to 2.71
- switch to pkg-readmes
- add rc scripts

pointers with README's and ok jasper@@
@
text
@d3 1
a3 1
# $OpenBSD$
d7 1
a7 1
daemon=${TRUEPREFIX}/bin/mimedefang
d12 4
@

