# cvs/ports/net/unifi/pkg/unifi.rc,v
head	1.6;
access;
symbols
	OPENBSD_6_1:1.4.0.6
	OPENBSD_6_1_BASE:1.4
	OPENBSD_6_0:1.4.0.4
	OPENBSD_6_0_BASE:1.4
	OPENBSD_5_9:1.4.0.2
	OPENBSD_5_9_BASE:1.4
	OPENBSD_5_8:1.3.0.6
	OPENBSD_5_8_BASE:1.3
	OPENBSD_5_7:1.3.0.2
	OPENBSD_5_7_BASE:1.3
	OPENBSD_5_6:1.2.0.2
	OPENBSD_5_6_BASE:1.2
	sthen_20140418:1.1.1.1
	sthen:1.1.1;
locks; strict;
comment	@ * @;


1.6
date	2017.05.27.19.57.08;	author ajacoutot;	state Exp;
branches;
next	1.5;
commitid	H08QxD8pqoInz6Rz;

1.5
date	2017.05.27.13.53.07;	author ajacoutot;	state Exp;
branches;
next	1.4;
commitid	SDpZdO8hfYIlMQ9j;

1.4
date	2016.02.12.13.46.45;	author sthen;	state Exp;
branches;
next	1.3;
commitid	6CWUg1LJQ9A1O3YS;

1.3
date	2014.10.03.11.53.23;	author jasper;	state Exp;
branches;
next	1.2;
commitid	YnCW657LSFnl923G;

1.2
date	2014.04.18.22.01.23;	author sthen;	state Exp;
branches;
next	1.1;

1.1
date	2014.04.18.21.49.44;	author sthen;	state Exp;
branches
	1.1.1.1;
next	;

1.1.1.1
date	2014.04.18.21.49.44;	author sthen;	state Exp;
branches;
next	;


desc
@@


1.6
log
@Bring rc_bg back.
@
text
@#!/bin/sh
#
# $OpenBSD: unifi.rc,v 1.5 2017/05/27 13:53:07 ajacoutot Exp $

daemon="${TRUEPREFIX}/share/unifi/lib/ace.jar"
daemon_user="_unifi"

. /etc/rc.d/rc.subr

java="$(${LOCALBASE}/bin/javaPathHelper -c unifi)"

pexp="${java}${daemon_flags:+ ${daemon_flags}} -jar ${daemon} start"
rc_bg=YES
rc_reload=NO

rc_start() {
	${rcexec} "${java} ${daemon_flags} -jar ${daemon} start"
}

rc_stop() {
	${java} -jar ${daemon} stop
}

rc_cmd $1
@


1.5
log
@Cope with recent rc.subr(8) change (removal of rc_bg etc.)
@
text
@d3 1
a3 1
# $OpenBSD: unifi.rc,v 1.4 2016/02/12 13:46:45 sthen Exp $
d11 1
d13 1
a13 1

@


1.4
log
@Fix unifi's rc script to use daemon_flags. Now you can do "rcctl set unifi
flags -Xmx512M -Xms256M" or similar if you want to impose some bounds on
java's memory usage.  OK/feedback aja@@
@
text
@d3 1
a3 1
# $OpenBSD: unifi.rc,v 1.3 2014/10/03 11:53:23 jasper Exp $
d16 1
a16 1
	${rcexec} "${java} ${daemon_flags} -jar ${daemon} start" &
@


1.3
log
@rcexec should only really be used in rc_start, unless there's a good reason to
use it elsewhere

prompted by aja@@
ok aja@@ sthen@@
@
text
@d3 1
a3 1
# $OpenBSD: unifi.rc,v 1.2 2014/04/18 22:01:23 sthen Exp $
d11 1
a11 1
pexp="${java} -jar ${daemon} start"
d16 1
a16 1
	${rcexec} "${java} -jar ${daemon} start" &
@


1.2
log
@remove a useless rc_bg, from aja@@
@
text
@d3 1
a3 1
# $OpenBSD: unifi.rc,v 1.1.1.1 2014/04/18 21:49:44 sthen Exp $
d20 1
a20 1
	${rcexec} "${java} -jar ${daemon} stop"
@


1.1
log
@Initial revision
@
text
@d3 1
a3 1
# $OpenBSD$
a12 1
rc_bg=YES
@


1.1.1.1
log
@Import ports/net/unifi, controller software for ubiquiti unifi access points.
This is v3.1.10 which is not a full release yet but working stably for me.
ok jca@@
@
text
@@
