# cvs/ports/sysutils/collectd/pkg/collectd.rc,v
head	1.2;
access;
symbols
	OPENBSD_6_0:1.2.0.6
	OPENBSD_6_0_BASE:1.2
	OPENBSD_5_9:1.2.0.2
	OPENBSD_5_9_BASE:1.2
	OPENBSD_5_8:1.2.0.4
	OPENBSD_5_8_BASE:1.2
	OPENBSD_5_7:1.1.0.2
	OPENBSD_5_7_BASE:1.1;
locks; strict;
comment	@ * @;


1.2
date	2015.07.23.07.33.17;	author landry;	state Exp;
branches;
next	1.1;
commitid	mFWJoll0XLJ4Rcdj;

1.1
date	2014.09.25.20.10.13;	author landry;	state Exp;
branches;
next	;
commitid	ZsuZqR6aoW7tRgVV;


desc
@@


1.2
log
@Update to collectd 5.5.0.

Enable aggregation, curl_json, fhcount, log_logstash, memcached,
threshold, write_log & write_riemann plugins.
Add dependencies on devel/libyajl for curl_json & log_logstash and
devel/protobuf-c for write_riemann.
Add a README explaining how to run the daemon as _collectd user in case
you dont need the dns and ping plugins (the only ones requiring root).
Lots of patches merged upstream.

From Pavel Korovin with small tweaks by me, tested here on ppc.
@
text
@#!/bin/sh
#
# $OpenBSD: collectd.rc,v 1.1 2014/09/25 20:10:13 landry Exp $

daemon="${TRUEPREFIX}/sbin/collectd"

# Ensure _collectd can read ${SYSCONFDIR}/collectd.conf
rc_pre() {
    chgrp _collectd /etc/collectd.conf
    chmod g+r /etc/collectd.conf
}

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
@


1.1
log
@Long awaited update to collectd 5.4.1.

- remove useless MESSAGE, provide an rc script instead
- enable pf, target_v5upgrade, tail_csv and write_graphite plugins
- backport fixes from FreeBSD and NetBSD to build with newer libstatgrab
  (upstream collectd github ticket #445)
- remove the check and inclusion of pcap-bpf.h, we dont have/need it
- users of the previous 4.x version should refer to
  https://collectd.org/wiki/index.php/V4_to_v5_migration_guide for the
upgrade path.

tested on ppc & amd64.
@
text
@d3 1
a3 1
# $OpenBSD$
d6 6
@

