# cvs/ports/lang/php/5.6/pkg/Attic/php_fpm.rc,v
head	1.3;
access;
symbols
	OPENBSD_5_8:1.2.0.4
	OPENBSD_5_8_BASE:1.2
	OPENBSD_5_7:1.1.0.2
	OPENBSD_5_7_BASE:1.1;
locks; strict;
comment	@ * @;


1.3
date	2015.09.18.09.18.23;	author sthen;	state dead;
branches;
next	1.2;
commitid	x1OqzrB6vMfHx6Ez;

1.2
date	2015.07.15.16.03.07;	author ajacoutot;	state Exp;
branches;
next	1.1;
commitid	bHx0hnbnAsyWwZso;

1.1
date	2014.11.15.13.43.02;	author robert;	state Exp;
branches;
next	;
commitid	qky5fmLhNW6GE6fd;


desc
@@


1.3
log
@de-conflict php-fpm packages, ok robert@@

Heads up! The rc script has changed name:

  sed -i s/php_fpm/php56_fpm/g /etc/rc.conf.local
@
text
@#!/bin/sh
#
# $OpenBSD: php_fpm.rc,v 1.2 2015/07/15 16:03:07 ajacoutot Exp $

daemon="${TRUEPREFIX}/sbin/php-fpm-${PV}"

. /etc/rc.d/rc.subr

pexp="php-fpm-${PV}: master process .*"
rc_reload=NO

rc_cmd $1
@


1.2
log
@Tweak pexp a bit.
@
text
@d3 1
a3 1
# $OpenBSD: php_fpm.rc,v 1.1 2014/11/15 13:43:02 robert Exp $
@


1.1
log
@add php-5.6; ok sthen@@
@
text
@d3 1
a3 1
# $OpenBSD: php_fpm.rc,v 1.1.1.1 2013/08/07 16:46:56 robert Exp $
d9 1
a9 1
pexp="php-fpm-${PV}: master.*"
@

