# cvs/ports/games/yquake2/pkg/q2ded.rc,v
head	1.3;
access;
symbols
	OPENBSD_6_2:1.3.0.2
	OPENBSD_6_2_BASE:1.3
	OPENBSD_6_1:1.1.0.4
	OPENBSD_6_1_BASE:1.1
	OPENBSD_6_0:1.1.0.2
	OPENBSD_6_0_BASE:1.1;
locks; strict;
comment	@ * @;


1.3
date	2017.05.27.19.57.06;	author ajacoutot;	state Exp;
branches;
next	1.2;
commitid	H08QxD8pqoInz6Rz;

1.2
date	2017.05.27.13.53.06;	author ajacoutot;	state Exp;
branches;
next	1.1;
commitid	SDpZdO8hfYIlMQ9j;

1.1
date	2016.06.28.20.50.22;	author awolk;	state Exp;
branches;
next	;
commitid	OCjzSGt5M9IIpgwf;


desc
@@


1.3
log
@Bring rc_bg back.
@
text
@#!/bin/sh
#
# $OpenBSD: q2ded.rc,v 1.2 2017/05/27 13:53:06 ajacoutot Exp $

daemon="${TRUEPREFIX}/bin/q2ded"
daemon_user="_yquake2"

. /etc/rc.d/rc.subr

pexp="${daemon}.*"
rc_bg=YES
rc_reload=NO

rc_cmd $1
@


1.2
log
@Cope with recent rc.subr(8) change (removal of rc_bg etc.)
@
text
@d3 1
a3 1
# $OpenBSD: q2ded.rc,v 1.1 2016/06/28 20:50:22 awolk Exp $
d11 1
a11 1

@


1.1
log
@Update games/yquake2 5.32 => 5.34

OK abieber@@

Changes in the port:
 - taking MAINTAINER, agreed with bmercer@@ who provided the initial port
 - added a q2ded rc script in a similar fashion as the games/ioquake3
   port does
 - the existing /usr/local/share/yquake2/q2ded binary from previous
   port did not work as it was searching for games.so from a relative
   path. I found a way to make both quake2 and q2ded work with an
   upstream supported MAKE variable to provide the system wide quake
   data directory. Thanks to this change we can drop the wrapper script
   and use the binaries directly.
 - adding a new _yquake2 user/group plus /var/yquake2 for the dedicated server
   binary
 - updated the port to version 5.34

Additional notes:
 - make port-lib-depends-check complains on:
   Extra:  openal.2 pthread.22

   ktrace tells me that the check is wrong. Openal and pthreads are both
   used.
 - upstream said in their changelog that they are changing their
   platform policy 'Switch from an arch whitelist to an "all archs are
   supported" approach.'. I saw in our tree a change by landry@@
   specifically limiting the supported platform to the upstream
   supported versions
   (http://cvsweb.openbsd.org/cgi-bin/cvsweb/ports/games/yquake2/Makefile?rev=1.2&content-type=text/x-cvsweb-markup).
   I am not removing this restriction until I get some feedback from
   people able to test other platforms than currently
   listed in ONLY_FOR_ARCHS.

Upstream changelog:

Quake II 5.33 to 5.34:
- Add support for stereo 3D (by Valery Guskov)
- Make gibt solid so they move with conveyor belts.
- Disable gl_ext_multitexturing by default.
- Switch from an arch whitelist to an "all archs are supported"
  approach.
- Add a new README.

Quake II 5.32 to 5.33:
- Add OGG volume slider in settings menu
- Fixed some bugs in volume settings
- Replaced HUD scale option in video menu with generic UI scale
- General UI upscaling improvements
- Better support for keyboards with AZERTY layout
@
text
@d3 1
a3 1
# $OpenBSD: $
a11 1
rc_bg=YES
@

