# cvs/ports/mail/spampd/pkg/spampd.rc,v
head	1.2;
access;
symbols
	OPENBSD_6_2:1.2.0.12
	OPENBSD_6_2_BASE:1.2
	OPENBSD_6_1:1.2.0.10
	OPENBSD_6_1_BASE:1.2
	OPENBSD_6_0:1.2.0.8
	OPENBSD_6_0_BASE:1.2
	OPENBSD_5_9:1.2.0.4
	OPENBSD_5_9_BASE:1.2
	OPENBSD_5_8:1.2.0.6
	OPENBSD_5_8_BASE:1.2
	OPENBSD_5_7:1.2.0.2
	OPENBSD_5_7_BASE:1.2
	OPENBSD_5_6:1.1.1.1.0.4
	OPENBSD_5_6_BASE:1.1.1.1
	OPENBSD_5_5:1.1.1.1.0.2
	OPENBSD_5_5_BASE:1.1.1.1
	ajacoutot_20140208:1.1.1.1
	ajacoutot:1.1.1;
locks; strict;
comment	@ * @;


1.2
date	2014.11.27.12.37.29;	author ajacoutot;	state Exp;
branches;
next	1.1;
commitid	jupCUYRTaGIqa8NY;

1.1
date	2014.02.08.18.58.22;	author ajacoutot;	state Exp;
branches
	1.1.1.1;
next	;

1.1.1.1
date	2014.02.08.18.58.22;	author ajacoutot;	state Exp;
branches;
next	;


desc
@@


1.2
log
@Bring some patches from mpaperno/spampd on github to fix some runtime issues
reported by Wesley Mouedine Assaby.
@
text
@#!/bin/sh
#
# $OpenBSD: spampd.rc,v 1.1.1.1 2014/02/08 18:58:22 ajacoutot Exp $

daemon="${TRUEPREFIX}/sbin/spampd"
daemon_flags="--a --aw --p=/var/spampd/spampd.pid"

. /etc/rc.d/rc.subr

pexp="/usr/bin/perl -T ${daemon}${daemon_flags:+ ${daemon_flags}}"

rc_cmd $1
@


1.1
log
@Initial revision
@
text
@d3 1
a3 1
# $OpenBSD$
d6 1
a6 1
daemon_flags="--tagall --aw"
@


1.1.1.1
log
@Import spampd-2.30.

spampd is a program used within an e-mail delivery system to scan messages
for possible Unsolicited Commercial E-mail (UCE, aka spam) content.
<...>

based on a submission by Henrik Friedrichsen (maintainer)
"get it in" espie@@
ok rpointel@@
@
text
@@
