# cvs/ports/net/xl2tpd/pkg/xl2tpd.rc,v
head	1.2;
access;
symbols
	OPENBSD_6_0:1.2.0.8
	OPENBSD_6_0_BASE:1.2
	OPENBSD_5_9:1.2.0.4
	OPENBSD_5_9_BASE:1.2
	OPENBSD_5_8:1.2.0.6
	OPENBSD_5_8_BASE:1.2
	OPENBSD_5_7:1.2.0.2
	OPENBSD_5_7_BASE:1.2
	OPENBSD_5_6:1.1.1.1.0.8
	OPENBSD_5_6_BASE:1.1.1.1
	OPENBSD_5_5:1.1.1.1.0.6
	OPENBSD_5_5_BASE:1.1.1.1
	OPENBSD_5_4:1.1.1.1.0.4
	OPENBSD_5_4_BASE:1.1.1.1
	OPENBSD_5_3:1.1.1.1.0.2
	OPENBSD_5_3_BASE:1.1.1.1
	sthen_20130115:1.1.1.1
	sthen:1.1.1;
locks; strict;
comment	@ * @;


1.2
date	2014.10.11.21.11.20;	author sthen;	state Exp;
branches;
next	1.1;
commitid	o8l9PCL4q7XU0XHl;

1.1
date	2013.01.15.23.02.19;	author sthen;	state Exp;
branches
	1.1.1.1;
next	;

1.1.1.1
date	2013.01.15.23.02.19;	author sthen;	state Exp;
branches;
next	;


desc
@@


1.2
log
@don't re-create the fifo if it still exists from a previous failed startup,
noted by landry
@
text
@#!/bin/sh
#
# $OpenBSD: xl2tpd.rc,v 1.1.1.1 2013/01/15 23:02:19 sthen Exp $

daemon="${TRUEPREFIX}/sbin/xl2tpd"
daemon_user="_xl2tpd"

. /etc/rc.d/rc.subr

rc_pre() {
	/usr/bin/install -d -o ${daemon_user} /var/run/xl2tpd
	[ -p /var/run/xl2tpd/l2tp-control ] ||
	    mkfifo /var/run/xl2tpd/l2tp-control
}

rc_cmd $1
@


1.1
log
@Initial revision
@
text
@d3 1
a3 1
# $OpenBSD: mrtg.rc,v 1.5 2011/12/17 22:01:13 sthen Exp $
d12 2
a13 1
	mkfifo /var/run/xl2tpd/l2tp-control
@


1.1.1.1
log
@import ports/net/xl2tpd, ok pirofti@@

xl2tpd is an implementation of the Layer 2 Tunnelling Protocol (RFC 2661),
which works by opening a pseudo-tty for communicating with pppd.

L2TP allows you to tunnel PPP over UDP. Some ISPs use L2TP to tunnel user
sessions from dial-in servers (modem banks, ADSL DSLAMs) to back-end PPP
servers. Another important application is Virtual Private Networks where
the IPsec protocol is used to secure the L2TP connection (L2TP/IPsec,
RFC 3193).
@
text
@@
