# cvs/ports/net/bird/pkg/bird.rc,v
head	1.5;
access;
symbols
	OPENBSD_6_1:1.5.0.8
	OPENBSD_6_1_BASE:1.5
	OPENBSD_6_0:1.5.0.6
	OPENBSD_6_0_BASE:1.5
	OPENBSD_5_9:1.5.0.2
	OPENBSD_5_9_BASE:1.5
	OPENBSD_5_8:1.5.0.4
	OPENBSD_5_8_BASE:1.5
	OPENBSD_5_7:1.4.0.4
	OPENBSD_5_7_BASE:1.4
	OPENBSD_5_6:1.4.0.2
	OPENBSD_5_6_BASE:1.4
	OPENBSD_5_5:1.2.0.6
	OPENBSD_5_5_BASE:1.2
	OPENBSD_5_4:1.2.0.4
	OPENBSD_5_4_BASE:1.2
	OPENBSD_5_3:1.2.0.2
	OPENBSD_5_3_BASE:1.2
	OPENBSD_5_2:1.1.0.6
	OPENBSD_5_2_BASE:1.1
	OPENBSD_5_1_BASE:1.1
	OPENBSD_5_1:1.1.0.4
	OPENBSD_5_0:1.1.0.2
	OPENBSD_5_0_BASE:1.1;
locks; strict;
comment	@ * @;


1.5
date	2015.07.15.14.35.48;	author ajacoutot;	state Exp;
branches;
next	1.4;
commitid	KYItW6R2sTqvHRL5;

1.4
date	2014.07.01.21.51.42;	author sthen;	state Exp;
branches;
next	1.3;
commitid	n8getCB40HHpwti4;

1.3
date	2014.07.01.20.38.44;	author sthen;	state Exp;
branches;
next	1.2;
commitid	2QZTmde4SPnd8qJ8;

1.2
date	2013.01.24.10.04.15;	author sthen;	state Exp;
branches
	1.2.6.1;
next	1.1;

1.1
date	2011.05.02.10.49.46;	author sthen;	state Exp;
branches;
next	;

1.2.6.1
date	2014.07.01.20.40.40;	author sthen;	state Exp;
branches;
next	;
commitid	AGvZZkTzeZ6DgMux;


desc
@@


1.5
log
@Adapt pexp after recent rc.subr change.
@
text
@#!/bin/sh
#
# $OpenBSD: bird.rc,v 1.4 2014/07/01 21:51:42 sthen Exp $

daemon="${TRUEPREFIX}/sbin/bird"

. /etc/rc.d/rc.subr

rc_stop() {
	${TRUEPREFIX}/sbin/birdc down
}

rc_cmd $1
@


1.4
log
@no need for the custom rc_check if pexp is put in the right place ;)
clue from rpe@@
@
text
@d3 1
a3 1
# $OpenBSD: bird.rc,v 1.3 2014/07/01 20:38:44 sthen Exp $
a7 2

pexp="${daemon}${daemon_flags:+ ${daemon_flags}}$"
@


1.3
log
@avoid rc.d/bird from matching bird6 with its pexp line, problem reported by
Job Snijders via jasper@@
@
text
@d3 1
a3 1
# $OpenBSD: bird.rc,v 1.2 2013/01/24 10:04:15 sthen Exp $
d9 1
a9 3
rc_check() {
        pgrep -q -f "^${pexp}$"
}
@


1.2
log
@Add comments about flavour handling to makefile.
Use birdc/birdc6 down for rc_stop.
@
text
@d3 1
a3 1
# $OpenBSD: bird.rc,v 1.1 2011/05/02 10:49:46 sthen Exp $
d8 4
@


1.2.6.1
log
@avoid rc.d/bird from matching bird6 with its pexp line, problem reported by
Job Snijders via jasper@@
@
text
@d3 1
a3 1
# $OpenBSD: bird.rc,v 1.3 2014/07/01 20:38:44 sthen Exp $
a7 4

rc_check() {
        pgrep -q -f "^${pexp}$"
}
@


1.1
log
@- update to 1.3.1 (bug fixes + adds 'med metric' for bgp)
- add rc scripts
@
text
@d3 1
a3 1
# $OpenBSD: mrtg.rc,v 1.4 2011/03/20 06:43:47 ajacoutot Exp $
d8 4
@

