# cvs/ports/mail/mimedefang/pkg/mimedefang.rc,v
head	1.2;
access;
symbols
	OPENBSD_6_0:1.2.0.14
	OPENBSD_6_0_BASE:1.2
	OPENBSD_5_9:1.2.0.10
	OPENBSD_5_9_BASE:1.2
	OPENBSD_5_8:1.2.0.12
	OPENBSD_5_8_BASE:1.2
	OPENBSD_5_7:1.2.0.8
	OPENBSD_5_7_BASE:1.2
	OPENBSD_5_6:1.2.0.6
	OPENBSD_5_6_BASE:1.2
	OPENBSD_5_5:1.2.0.4
	OPENBSD_5_5_BASE:1.2
	OPENBSD_5_4:1.2.0.2
	OPENBSD_5_4_BASE:1.2
	OPENBSD_5_3:1.1.0.2
	OPENBSD_5_3_BASE:1.1;
locks; strict;
comment	@ * @;


1.2
date	2013.03.28.14.27.19;	author ajacoutot;	state Exp;
branches;
next	1.1;

1.1
date	2012.12.05.18.51.56;	author ajacoutot;	state Exp;
branches;
next	;


desc
@@


1.2
log
@Enable queuing by default (up to 10 msgs).
Rework and document the way mimedefang+spamassassin work wrt bayes.
Add a note on using mfs(8) for the spool directory.
Take maintainer.

ok okan@@ (former maintainer)
@
text
@#!/bin/sh
#
# $OpenBSD: mimedefang.rc,v 1.1 2012/12/05 18:51:56 ajacoutot Exp $

daemon="${TRUEPREFIX}/bin/mimedefang"
daemon_flags="-Gqr -U _mdefang -P /var/spool/MIMEDefang/mimedefang.pid -p /var/spool/MIMEDefang/mimedefang.sock -m /var/spool/MIMEDefang/mimedefang-multiplexor.sock"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
@


1.1
log
@Add `-r' to mimedefang flags by default.
Add an example for filter_relay() that works out of the box.
Rename the mimedefang_daemon rc.d(8) script to mimedefang as it's
supposed to be...

ok okan@@ (maintainer)
@
text
@d3 1
a3 1
# $OpenBSD: mimedefang_daemon.rc,v 1.7 2012/09/18 15:34:39 ajacoutot Exp $
d6 1
a6 1
daemon_flags="-Gr -U _mdefang -P /var/spool/MIMEDefang/mimedefang.pid -p /var/spool/MIMEDefang/mimedefang.sock -m /var/spool/MIMEDefang/mimedefang-multiplexor.sock"
@

