# cvs/ports/sysutils/logstash/forwarder/pkg/Attic/logstash_forwarder.rc,v
head	1.3;
access;
symbols
	OPENBSD_5_9:1.2.0.2
	OPENBSD_5_9_BASE:1.2
	OPENBSD_5_8:1.2.0.4
	OPENBSD_5_8_BASE:1.2
	OPENBSD_5_7:1.1.1.1.0.2
	OPENBSD_5_7_BASE:1.1.1.1
	jasper_20141612:1.1.1.1
	jasper:1.1.1;
locks; strict;
comment	@ * @;


1.3
date	2016.04.26.13.11.20;	author jasper;	state dead;
branches;
next	1.2;
commitid	arm9BxtVTyUl3xtv;

1.2
date	2015.05.15.07.37.14;	author jasper;	state Exp;
branches;
next	1.1;
commitid	gwLRAJ7a60eNtifR;

1.1
date	2014.12.16.15.21.06;	author jasper;	state Exp;
branches
	1.1.1.1;
next	;
commitid	S0ETu7llOHsw6Hm1;

1.1.1.1
date	2014.12.16.15.21.06;	author jasper;	state Exp;
branches;
next	;
commitid	S0ETu7llOHsw6Hm1;


desc
@@


1.3
log
@remove logstash-forwarder, it's no longer developed upstream. filebeat
should be used instead. for more information on migrating: https://www.elastic.co/guide/en/beats/filebeat/current/migrating-from-logstash-forwarder.html
@
text
@#!/bin/sh
#
# $OpenBSD: logstash_forwarder.rc,v 1.2 2015/05/15 07:37:14 jasper Exp $

daemon="${PREFIX}/sbin/logstash-forwarder"
daemon_flags="-config ${SYSCONFDIR}/logstash-forwarder/logstash-forwarder.conf"

. /etc/rc.d/rc.subr

rc_bg=YES
rc_reload=NO

rc_cmd $1
@


1.2
log
@- update to logstash-forwarder 0.4.0
* release notes at https://www.elastic.co/blog/logstash-forwarder-0-4-0-released
- sync default config/readme with logstash
@
text
@d3 1
a3 1
# $OpenBSD: logstash_forwarder.rc,v 1.1.1.1 2014/12/16 15:21:06 jasper Exp $
@


1.1
log
@Initial revision
@
text
@d3 1
a3 1
# $OpenBSD$
d6 1
a6 1
daemon_flags="-config ${SYSCONFDIR}/logstash-forwarder/logstash-forwarder.conf -v"
@


1.1.1.1
log
@import logstash-forwarder 0.3.1.1 (0.3.1 + git head)

Tool to collect logs locally in preparation for processing elsewhere with
logstash.

ok aja@@
@
text
@@
