# cvs/ports/sysutils/beats/heartbeat/pkg/heartbeat.rc,v
head     1.1;
branch   1.1.1;
access   ;
symbols  jasper_20172904:1.1.1.1 jasper:1.1.1;
locks    ; strict;
comment  @ * @;


1.1
date     2017.04.29.14.16.10;  author jasper;  state Exp;
branches 1.1.1.1;
next     ;
commitid        nftAqjuXKadn7ufN;

1.1.1.1
date     2017.04.29.14.16.10;  author jasper;  state Exp;
branches ;
next     ;
commitid        nftAqjuXKadn7ufN;


desc
@@



1.1
log
@Initial revision
@
text
@#!/bin/sh
#
# $OpenBSD: metricbeat.rc,v 1.1.1.1 2017/04/27 09:17:26 jasper Exp $

daemon="${PREFIX}/bin/heartbeat"
daemon_flags="-path.config ${SYSCONFDIR}/heartbeat/"

. /etc/rc.d/rc.subr

rc_bg=YES
rc_reload=NO

rc_cmd $1
@


1.1.1.1
log
@import heartbeat-5.2.1

Monitor services for their availability with active probing. Given a
list of URLs, Heartbeat asks the simple question: Are you alive?
Heartbeat ships this information and response time to the rest of the
Elastic Stack for further analysis.

ok aja@@
@
text
@@
