# cvs/ports/net/ejabberd/pkg/ejabberd.rc,v
head	1.2;
access;
symbols
	OPENBSD_6_2:1.2.0.12
	OPENBSD_6_2_BASE:1.2
	OPENBSD_6_1:1.2.0.10
	OPENBSD_6_1_BASE:1.2
	OPENBSD_6_0:1.2.0.8
	OPENBSD_6_0_BASE:1.2
	OPENBSD_5_9:1.2.0.4
	OPENBSD_5_9_BASE:1.2
	OPENBSD_5_8:1.2.0.6
	OPENBSD_5_8_BASE:1.2
	OPENBSD_5_7:1.2.0.2
	OPENBSD_5_7_BASE:1.2;
locks; strict;
comment	@ * @;


1.2
date	2014.12.10.15.30.19;	author ajacoutot;	state Exp;
branches;
next	1.1;
commitid	Sej0DMbFFetOhScw;

1.1
date	2014.12.10.15.25.05;	author sebastia;	state Exp;
branches;
next	;
commitid	YLBbSWvNDYP0RIRk;


desc
@@


1.2
log
@Cosmetic: follow rc.template.
@
text
@#!/bin/sh
#
# $OpenBSD: ejabberd.rc,v 1.1 2014/12/10 15:25:05 sebastia Exp $

daemon="${TRUEPREFIX}/sbin/ejabberdctl"

. /etc/rc.d/rc.subr

rc_reload=NO
rc_usercheck=NO

rc_start() {
	${rcexec} "${daemon} ${daemon_flags} start"
}

rc_stop() {
	${daemon} ${daemon_flags} stop
}

rc_check() {
	${daemon} ${daemon_flags} status
}

rc_cmd $1
@


1.1
log
@Add rcscript for ejabberd, update the README based on
feedback from zhuk@@, feedback and help with the rcscript
from aja@@

Remove MAINTAINER as per his request

OK ajacoutot@@
@
text
@d3 1
a3 1
# $OpenBSD$
d8 2
a22 2

rc_reload=NO
@

