# cvs/ports/net/pdnsd/pkg/pdnsd.rc,v
head	1.2;
access;
symbols
	OPENBSD_6_1:1.2.0.22
	OPENBSD_6_1_BASE:1.2
	OPENBSD_6_0:1.2.0.20
	OPENBSD_6_0_BASE:1.2
	OPENBSD_5_9:1.2.0.16
	OPENBSD_5_9_BASE:1.2
	OPENBSD_5_8:1.2.0.18
	OPENBSD_5_8_BASE:1.2
	OPENBSD_5_7:1.2.0.14
	OPENBSD_5_7_BASE:1.2
	OPENBSD_5_6:1.2.0.12
	OPENBSD_5_6_BASE:1.2
	OPENBSD_5_5:1.2.0.10
	OPENBSD_5_5_BASE:1.2
	OPENBSD_5_4:1.2.0.8
	OPENBSD_5_4_BASE:1.2
	OPENBSD_5_3:1.2.0.6
	OPENBSD_5_3_BASE:1.2
	OPENBSD_5_2:1.2.0.4
	OPENBSD_5_2_BASE:1.2
	OPENBSD_5_1_BASE:1.2
	OPENBSD_5_1:1.2.0.2;
locks; strict;
comment	@ * @;


1.2
date	2011.12.17.22.01.13;	author sthen;	state Exp;
branches;
next	1.1;

1.1
date	2011.10.14.07.00.51;	author ajacoutot;	state Exp;
branches;
next	;


desc
@@


1.2
log
@Where a program with an rc script _requires_ a certain command line
parameter to daemonize, move the parameter from daemon to daemon_flags,
so that the user cannot inadvertently prevent it from daemonizing by
adjusting the flags.

Discussed with ajacoutot and schwarze, this method was suggested
by schwarze@@ as a simpler alternative to my diff. ok aja@@
@
text
@#!/bin/sh
#
# $OpenBSD: pdnsd.rc,v 1.1 2011/10/14 07:00:51 ajacoutot Exp $

daemon="${TRUEPREFIX}/sbin/pdnsd -d"

. /etc/rc.d/rc.subr

rc_reload() {
	${TRUEPREFIX}/sbin/pdnsd-ctl config
}

rc_cmd $1
@


1.1
log
@Add an rc script, from frantisek holop.
Tweak the port a bit while here.
@
text
@d3 1
a3 1
# $OpenBSD$
d5 1
a5 2
daemon="${TRUEPREFIX}/sbin/pdnsd"
daemon_flags="-d"
@

