# cvs/ports/www/apache-httpd-openbsd/pkg/Attic/apache.rc,v
head	1.3;
access;
symbols
	OPENBSD_5_9:1.2.0.2
	OPENBSD_5_9_BASE:1.2
	OPENBSD_5_8:1.2.0.4
	OPENBSD_5_8_BASE:1.2
	OPENBSD_5_7:1.1.0.4
	OPENBSD_5_7_BASE:1.1
	OPENBSD_5_6:1.1.0.2
	OPENBSD_5_6_BASE:1.1;
locks; strict;
comment	@ * @;


1.3
date	2016.04.28.18.19.47;	author sthen;	state dead;
branches;
next	1.2;
commitid	IzB6dS8YEKbH0xwR;

1.2
date	2015.07.17.23.58.25;	author ajacoutot;	state Exp;
branches;
next	1.1;
commitid	ueY4CFFCF5BT06rk;

1.1
date	2014.07.27.10.51.37;	author kili;	state Exp;
branches;
next	;
commitid	gYt5lPsCZtEp8W4S;


desc
@@


1.3
log
@remove, nothing uses it any more
@
text
@#!/bin/sh
#
# $OpenBSD: apache.rc,v 1.2 2015/07/17 23:58:25 ajacoutot Exp $

daemon="${TRUEPREFIX}/sbin/httpd"

. /etc/rc.d/rc.subr

pexp="httpd: parent.*"

rc_cmd $1
@


1.2
log
@RUN_DEPENDS on www/apache-httpd,-common to get the icons.
@
text
@d3 1
a3 1
# $OpenBSD: apache.rc,v 1.1 2014/07/27 10:51:37 kili Exp $
@


1.1
log
@Rename rc scripts of apache-httpd to apache resp. apache2 to get rid
of a conflict with base httpd.

ok sthen@@
@
text
@d3 1
a3 1
# $OpenBSD$
a9 4

rc_pre() {
	rm -f /var/www/logs/{ssl_mutex,httpd.lock,accept.lock}.*
}
@

