# cvs/ports/www/apache-httpd-openbsd/pkg/Attic/httpd.rc,v
head	1.2;
access;
symbols
	sthen_20140313:1.1.1.1
	sthen:1.1.1;
locks; strict;
comment	@ * @;


1.2
date	2014.07.27.10.51.37;	author kili;	state dead;
branches;
next	1.1;
commitid	gYt5lPsCZtEp8W4S;

1.1
date	2014.03.13.17.28.04;	author sthen;	state Exp;
branches
	1.1.1.1;
next	;

1.1.1.1
date	2014.03.13.17.28.04;	author sthen;	state Exp;
branches;
next	;


desc
@@


1.2
log
@Rename rc scripts of apache-httpd to apache resp. apache2 to get rid
of a conflict with base httpd.

ok sthen@@
@
text
@#!/bin/sh
#
# $OpenBSD: httpd.rc,v 1.1 2014/03/13 17:28:04 sthen Exp $

daemon="${TRUEPREFIX}/sbin/httpd"

. /etc/rc.d/rc.subr

pexp="httpd: parent.*"

rc_pre() {
	rm -f /var/www/logs/{ssl_mutex,httpd.lock,accept.lock}.*
}

rc_cmd $1
@


1.1
log
@Initial revision
@
text
@d3 1
a3 1
# $OpenBSD: httpd,v 1.1 2011/07/06 18:55:36 robert Exp $
@


1.1.1.1
log
@Import a snapshot of the modified Apache httpd (previously in base) to ports.
Mostly worked on by florian@@, port tweaks by me.
@
text
@@
