# cvs/ports/net/geomyidae/pkg/geomyidae.rc,v
head	1.2;
access;
symbols
	OPENBSD_6_1:1.2.0.8
	OPENBSD_6_1_BASE:1.2
	OPENBSD_6_0:1.2.0.6
	OPENBSD_6_0_BASE:1.2
	OPENBSD_5_9:1.2.0.2
	OPENBSD_5_9_BASE:1.2
	OPENBSD_5_8:1.2.0.4
	OPENBSD_5_8_BASE:1.2
	jung-20150622:1.1.1.1
	jung:1.1.1;
locks; strict;
comment	@ * @;


1.2
date	2015.07.15.19.27.52;	author ajacoutot;	state Exp;
branches;
next	1.1;
commitid	vKlnNBj0UJE6itee;

1.1
date	2015.06.22.15.53.45;	author jung;	state Exp;
branches
	1.1.1.1;
next	;
commitid	PI2jkSusjTWJBWyH;

1.1.1.1
date	2015.06.22.15.53.45;	author jung;	state Exp;
branches;
next	;
commitid	PI2jkSusjTWJBWyH;


desc
@@


1.2
log
@Fix pexp.
@
text
@#!/bin/sh
#
# $OpenBSD: geomyidae.rc,v 1.1.1.1 2015/06/22 15:53:45 jung Exp $

daemon="${TRUEPREFIX}/bin/geomyidae"
daemon_flags="-4 -l /var/log/geomyidae.log -u _geomyidae -g _geomyidae"

. /etc/rc.d/rc.subr

pexp="${daemon} .*"

rc_cmd $1
@


1.1
log
@Initial revision
@
text
@d3 1
a3 1
# $OpenBSD: rc.template,v 1.6 2014/08/14 19:54:10 ajacoutot Exp $
d10 1
a10 1
pexp="${daemon}.*"
@


1.1.1.1
log
@import net/geomyidae v0.26.3
 
Geomyidae is a daemon for serving the Gopher protocol.
It is based on gopherd for Plan 9.

ok jca@@

@
text
@@
